# poky/meta/recipes-extended/sysklogd/files/klogd.service
[Unit]
Description=Kernel Logging Service
After=syslogd.service

[Service]
Type=forking
ExecStart=/sbin/klogd
PIDFile=/var/run/klogd.pid
StandardOutput=null
Restart=on-failure

[Install]
WantedBy=multi-user.target
