# meta-phosphor/common/recipes-phosphor/settings/phosphor-settings/org.openbmc.settings.Host.service
[Unit]
Description=Phosphor Settings Manager
Before=mapper-wait@-org-openbmc-settings-host0.service
Wants=obmc-mapper.target
After=obmc-mapper.target

[Service]
ExecStart={sbindir}/settings_manager.py
Environment="PYTHONUNBUFFERED=1"
Restart=always
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
