# meta-fii/meta-kudo/recipes-kudo/kudo-sys-utility/kudo-boot/host-powerctrl.service
[Unit]
Description = Ampere Host Power Control Sequence

[Service]
Restart=no
RemainAfterExit=true
Type=oneshot
ExecStart=/usr/bin/init_once.sh
StandardOutput=syslog

[Install]
WantedBy=sysinit.target
