# meta-fii/meta-kudo/recipes-kudo/kudo-sys-utility/kudo-boot/host-powerctrl.service
[Unit]
Description = Host Power Control Sequence
Wants=ssifbridge.service
After=ssifbridge.service
Wants=kudo-bios-verify.service
After=kudo-bios-verify.service

[Service]
Restart=no
RemainAfterExit=true
Type=oneshot
ExecStart=/usr/libexec/kudo-boot/init_once.sh
StandardOutput=syslog

[Install]
WantedBy=sysinit.target
