# meta-phosphor/recipes-core/systemd/obmc-targets/obmc-power-on@.target
[Unit]
Description=Power%i (On)
Wants={SYSTEMD_DEFAULT_TARGET}
After=obmc-power-start@%i.target
Conflicts=obmc-chassis-poweroff@%i.target
RefuseManualStart=yes
RefuseManualStop=yes
