# meta-phosphor/recipes-core/systemd/obmc-targets/obmc-host-startmin@.target
[Unit]
Description=Start Host%i Minimum
Wants=multi-user.target
After=multi-user.target
Conflicts=obmc-host-stop@%i.target
RefuseManualStop=yes
OnFailure=obmc-host-quiesce@%i.target
OnFailureJobMode=flush
