# meta-phosphor/recipes-core/systemd/obmc-targets/obmc-fans-ready.target
[Unit]
Description=Fans Ready
RefuseManualStart=yes
RefuseManualStop=yes
