# meta-phosphor/recipes-core/systemd/obmc-targets/obmc-power-reset-on@.target
[Unit]
Description=Chassis%i power on after reset
Conflicts=obmc-chassis-poweroff@%i.target
RefuseManualStart=yes
RefuseManualStop=yes
