# meta-phosphor/recipes-core/systemd/obmc-targets/obmc-host-quiesce@.target
[Unit]
Description=Quiesce Target
Wants=multi-user.target
After=multi-user.target
RefuseManualStop=yes
Conflicts=obmc-chassis-poweroff@%i.target
