# meta-phosphor/recipes-core/systemd/obmc-targets/obmc-host-crash@.target
[Unit]
Description=Host instance %i crashed
Conflicts=obmc-chassis-poweroff@%i.target
RefuseManualStop=yes
