# meta-phosphor/recipes-core/systemd/obmc-targets/obmc-host-start-pre@.target
[Unit]
Description=Start Host%i (Pre)
Wants=multi-user.target
After=multi-user.target
Conflicts=obmc-host-stop@%i.target
RefuseManualStart=yes
RefuseManualStop=yes
