# meta-phosphor/recipes-core/systemd/obmc-targets/obmc-chassis-poweroff@.target
[Unit]
Description=Chassis%i (Power Off)
Wants={SYSTEMD_DEFAULT_TARGET}
After={SYSTEMD_DEFAULT_TARGET}
Wants=mapper-wait@-xyz-openbmc_project-state-chassis%i.service
After=mapper-wait@-xyz-openbmc_project-state-chassis%i.service
Conflicts=obmc-chassis-poweron@%i.target
Conflicts=obmc-chassis-powerreset@%i.target
RefuseManualStop=yes
