# meta-phosphor/recipes-core/systemd/obmc-targets/obmc-host-quiesce@.target
[Unit]
Description=Quiesce Target
Wants={SYSTEMD_DEFAULT_TARGET}
After={SYSTEMD_DEFAULT_TARGET}
RefuseManualStop=yes
Conflicts=obmc-chassis-poweroff@%i.target
