# meta-phosphor/recipes-core/systemd/obmc-targets/obmc-power-start-pre@.target
[Unit]
Description=Power%i On (Pre)
Wants=multi-user.target
After=multi-user.target
Conflicts=obmc-chassis-poweroff@%i.target
RefuseManualStart=yes
RefuseManualStop=yes
