# meta-phosphor/recipes-core/systemd/obmc-targets/obmc-power-stop-pre@.target
[Unit]
Description=Power%i Off (Pre)
Wants=multi-user.target
After=multi-user.target
Wants=obmc-host-stopped@%i.target
After=obmc-host-stopped@%i.target
Conflicts=obmc-chassis-poweron@%i.target
RefuseManualStart=yes
RefuseManualStop=yes
