# meta-phosphor/recipes-core/systemd/obmc-targets/obmc-host-stop@.target
[Unit]
Description=Stop Host%i
Wants=multi-user.target
After=multi-user.target
Wants=mapper-wait@-org-openbmc-control-chassis%i.service
After=mapper-wait@-org-openbmc-control-chassis%i.service
After=obmc-chassis-poweroff@%i.target
Conflicts=obmc-host-startmin@%i.target
RefuseManualStop=yes
OnFailure=obmc-chassis-poweroff@%i.target
OnFailureJobMode=flush
