# meta-phosphor/recipes-core/systemd/obmc-targets/obmc-host-stop-pre@.target
[Unit]
Description=Stop Host%i (Pre)
Wants=multi-user.target
After=multi-user.target
Before=obmc-host-stopping@%i.target
Conflicts=obmc-host-startmin@%i.target
RefuseManualStart=yes
RefuseManualStop=yes
