# meta-phosphor/recipes-core/systemd/obmc-targets/obmc-host-stopped@.target
[Unit]
Description=Host%i (Stopped)
Wants=multi-user.target
Wants=obmc-host-stopping@%i.target
After=obmc-host-stopping@%i.target
Before=obmc-power-stop-pre@%i.target
Conflicts=obmc-host-startmin@%i.target
RefuseManualStart=yes
RefuseManualStop=yes
