# meta-phosphor/recipes-core/systemd/obmc-targets/obmc-host-shutdown@.target
[Unit]
Description=Power%i Host Off
Wants=multi-user.target
RefuseManualStop=yes
OnFailure=obmc-chassis-poweroff@%i.target
OnFailureJobMode=flush
