# meta-phosphor/recipes-core/systemd/obmc-targets/obmc-host-shutdown@.target
[Unit]
Description=Power%i Host Off
Wants={SYSTEMD_DEFAULT_TARGET}
RefuseManualStop=yes
OnFailure=obmc-chassis-poweroff@%i.target
OnFailureJobMode=flush
