# meta-phosphor/recipes-core/systemd/obmc-targets/obmc-host-reset@.target
[Unit]
Description=Host%i (Reset Check)
Conflicts=obmc-host-stop@%i.target
RefuseManualStop=yes

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
