# meta-phosphor/recipes-core/systemd/obmc-targets/obmc-power-start-pre@.target
[Unit]
Description=Power%i On (Pre)
Wants={SYSTEMD_DEFAULT_TARGET}
After={SYSTEMD_DEFAULT_TARGET}
Conflicts=obmc-chassis-poweroff@%i.target
RefuseManualStart=yes
RefuseManualStop=yes
