# meta-phosphor/recipes-core/systemd/obmc-targets/obmc-host-started@.target
[Unit]
Description=Host%i (Started)
Wants=multi-user.target
After=obmc-host-starting@%i.target
Conflicts=obmc-host-stop@%i.target
RefuseManualStart=yes
RefuseManualStop=yes
