# meta-phosphor/recipes-phosphor/flash/noverify-bmc-update/noverify-bmc-verify.service
[Unit]
Description=Don't verify the image, just stage it for update

[Service]
Type=oneshot
ExecStart=/bin/mv /run/initramfs/bmc-image /run/initramfs/image-bmc
