# meta-phosphor/common/recipes-phosphor/host-ipmid/host-ipmi-bt/host-ipmi-bt.service
[Unit]
Description=Phosphor OpenBMC BT to DBUS 

[Service]
Restart=always
ExecStart=/usr/sbin/host-ipmi-bt

[Install]
WantedBy=multi-user.target
