# meta-phosphor/common/recipes-phosphor/obmc-phosphor-chassis/obmc-watchdog/org.openbmc.watchdog.Host@.service
[Unit]
Description=Phosphor Watchdog%i

[Service]
Restart=always
ExecStart={sbindir}/host_watchdog.exe
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
