# meta-security/meta-tpm/recipes-tpm/trousers/files/tcsd.service
[Unit]
Description=TCG Core Services Daemon
After=syslog.target

[Service]
Type=forking
ExecStart=@SBINDIR@/tcsd

[Install]
WantedBy=multi-user.target
