# meta-openbmc-machines/meta-openpower/meta-ingrasys/meta-zaius/recipes-phosphor/chassis/vcs-control/vcs-off@.service
[Unit]
Description=Turn VCS rails off
Wants=obmc-power-stop-pre@%i.target
Before=obmc-power-stop-pre@%i.target
Conflicts=obmc-host-start@%i.target

[Service]
ExecStart={bindir}/zaius_vcs.sh off
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=obmc-host-stop@%i.target
