# yocto-poky/meta/recipes-graphics/x11-common/xserver-nodm-init/xserver-nodm.service
[Unit]
Description=Xserver startup without a display manager

[Service]
User=root
EnvironmentFile=/etc/default/xserver-nodm
ExecStart=/etc/X11/Xserver

[Install]
Alias=display-manager.service
