# meta-fii/meta-kudo/recipes-kudo/kudo-fw-utility/kudo-bios-update/phosphor-ipmi-flash-bios-update.service
[Unit]
Description=Phosphor-ipmi-flash update BIOS firmware service

[Service]
ExecStart=/usr/sbin/kudo-fw.sh bios /run/initramfs/bios-image
Type=oneshot
