# meta-nuvoton/recipes-nuvoton/loadmcu/files/mcu-version@.service
[Unit]
Description=NPCM7xx MCU F/W version service
Requires=xyz.openbmc_project.Software.BMC.Updater.service
Before=xyz.openbmc_project.Software.BMC.Updater.service

[Service]
Restart=no
Type=oneshot
ExecStart=/usr/bin/mcu-version.sh %i
SyslogIdentifier=mcu-version.sh

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
