# meta-nuvoton/recipes-nuvoton/loadmcu/files/mcu-update.service
[Unit]
Description=NPCM7xx MCU F/W update service

[Service]
ExecStart=/usr/bin/loadmcu -d /dev/mcu0 -s /tmp/image-mcu
Type=oneshot
ExecStopPost=/bin/systemctl start mcu-version@* --all
