# meta-openembedded/meta-networking/recipes-protocols/freediameter/files/freediameter.service
[Unit]
Description=freediameter daemon
After=network.target

[Service]
Type=simple
PIDFile=/run/freediameter.pid
ExecStart=@BINDIR@/freeDiameterd

[Install]
WantedBy=multi-user.target
