# meta-quanta/meta-gsj/recipes-gsj/detect-fan-fail/files/detect-fan-fail.service
[Unit]
Description=Detect Fan Fail Manager
After=phosphor-pid-control.service

[Service]
ExecStart=/usr/bin/detect-fan-fail.sh
Restart=always
StandardOutput=syslog

[Install]
WantedBy=multi-user.target
