# meta-ampere/meta-common/recipes-ac01/host/ampere-hostctrl/ampere-host-force-reset.service
[Unit]
Description=Ampere Computing LLC force reset host
Conflicts=obmc-chassis-poweroff@0.target
OnFailure=obmc-chassis-powerreset@0.target

[Service]
RemainAfterExit=no
Type=oneshot
ExecStart=/usr/bin/env ampere_power_util.sh mb force_reset
SyslogIdentifier=ampere_power_util.sh

