# meta-google/recipes-google/networking/gbmc-nic-config/gbmc-nic-dhcrelay.service.in
[Unit]
Description=gBMC DHCP Relay Agent Daemon
After=network.target
StartLimitIntervalSec=10
StartLimitBurst=3

[Service]
Restart=always
RestartSec=5
ExecStart=/usr/sbin/dhcrelay -d --no-pid -l gbmcnicdhcp @UPPERS@

[Install]
WantedBy=multi-user.target
