# meta-google/recipes-google/networking/gbmc-nic-config/gbmc-nic-ra@.service
[Unit]
Description=gBMC %i RA Discovery
After=network.target
StartLimitIntervalSec=10
StartLimitBurst=3

[Service]
Restart=always
RestartSec=5
ExecStart=/usr/libexec/gbmc-nic-ra.sh %i

[Install]
WantedBy=multi-user.target
