# meta-google/recipes-google/networking/files/gbmc-mac-config.service
[Unit]
Before=systemd-networkd.service

[Service]
Restart=on-failure
Type=oneshot
ExecStart=/usr/libexec/gbmc-mac-config.sh

[Install]
WantedBy=multi-user.target
