# meta-google/recipes-google/networking/files/gbmc-ip-monitor.service
[Unit]
Before=systemd-networkd.service

[Service]
Type=notify
ExecStart=/usr/libexec/gbmc-ip-monitor.sh

[Install]
WantedBy=multi-user.target
