# meta-google/recipes-google/networking/files/sslh.socket
[Unit]
Before=sslh.service

[Socket]
ListenStream=3967

[Install]
WantedBy=sockets.target
