# meta-phosphor/common/recipes-phosphor/rest-dbus/files/rest-dbus.service
[Unit]
Description=DBUS introspecting REST server.

[Service]
Restart=always
ExecStart={sbindir}/rest-dbus
Environment="PYTHONUNBUFFERED=1"
