# meta-phosphor/common/recipes-phosphor/rest-dbus/files/rest-dbus.socket
[Unit]
Description=DBUS introspecting REST server socket.

[Socket]
ListenStream=3000

[Install]
WantedBy=sockets.target
