# meta-phosphor/common/recipes-phosphor/ipmi/phosphor-ipmi-fru/obmc-read-eeprom@.service
[Unit]
Description=Read %I EEPROM
Wants=mapper-wait@-xyz-openbmc_project-inventory.service
After=mapper-wait@-xyz-openbmc_project-inventory.service

[Service]
Restart=no
Type=oneshot
EnvironmentFile={envfiledir}/obmc/eeproms/%I
ExecStart=/usr/bin/env phosphor-read-eeprom --eeprom $SYSFS_PATH --fruid $FRUID
SyslogIdentifier=phosphor-read-eeprom

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
