# meta-facebook/recipes-phosphor/sensors/dbus-sensors/xyz.openbmc_project.leakdetector.critical.assert@.service
[Unit]
Description=Leak detection critical assert handler for %i

[Service]
Type=oneshot
ExecStart=/usr/libexec/dbus-sensors/critical-leak-assert-handler %i
SyslogIdentifier=%i-critical-leak-assert-handler
