# meta-facebook/meta-bletchley/recipes-bletchley/fw-update/fw-update/flash-usb-controller@.service
[Unit]
Description=Flash usb controller image on SLED %i

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/sbin/usb-controller-update %i
