# upstream-layers/openembedded-core/meta/recipes-graphics/wayland/weston-init/weston.socket
[Unit]
Description=Weston socket
RequiresMountsFor=@runtimedir@

[Socket]
ListenStream=@runtimedir@/wayland-0
SocketMode=0775
SocketUser=@WESTON_USER@
SocketGroup=wayland
RemoveOnStop=yes

[Install]
WantedBy=sockets.target

