# meta-openembedded/meta-webserver/recipes-webadmin/netdata/netdata/netdata.service
# SPDX-License-Identifier: GPL-3.0-or-later
[Unit]
Description=Real time performance monitoring
RequiresMountsFor=/var
After=network.target

[Service]
Type=simple
ExecStartPre=/bin/mkdir -p /var/log/netdata
ExecStartPre=/bin/chown -R netdata:netdata /var/log/netdata
ExecStart=/usr/sbin/netdata -D -u netdata


[Install]
WantedBy=multi-user.target
