# meta-facebook/meta-tiogapass/recipes-extended/rsyslog/rsyslog/rotate-event-logs.service
[Unit]
Description=Rotates the event logs

[Service]
Type=oneshot
ExecStart=/usr/sbin/logrotate /etc/logrotate.d/logrotate.rsyslog

[Install]
WantedBy=multi-user.target
