# meta-facebook/meta-tiogapass/recipes-extended/rsyslog/rsyslog/rotate-event-logs.timer
[Unit]
Description=Run rotate-event-logs.service every minute

[Timer]
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target
