# meta-ampere/meta-jade/recipes-phosphor/gpio/tempevent-monitor/ampere_hightemp_stop.service
[Unit]
Description=Ampere Computing LLC Host hightemp event stop service
Conflicts=obmc-host-stop@%i.target

[Service]
Restart=no
ExecStart=/usr/bin/tempevent_log Hightemp Stop
