# meta-ibm/meta-system1/recipes-phosphor/flash/phosphor-software-manager/obmc-flash-host-bios@.service
[Unit]
Description=Flash Host Bios firmware image %I to Host
Wants=pch-standby.service
After=pch-standby.service

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/libexec/bios-update.sh /tmp/images/%i
