# meta-ibm/recipes-phosphor/leds/phosphor-led-manager/obmc-led-create-virtual-leds@.service
[Unit]
Description=Create the virtual physical LEDs

[Service]
Restart=always
ExecStart=/usr/sbin/phosphor-ledcontroller -p %f
SyslogIdentifier=obmc-led-create-virtual-leds
