# meta-phosphor/common/recipes-phosphor/obmc-phosphor-policy/files/obmc-phosphor-policyd.service
[Unit]
Description=Phosphor OpenBMC policy management daemon

[Service]
Restart=always
ExecStart=/usr/sbin/obmc-phosphor-policyd

[Install]
WantedBy=multi-user.target
