# meta-openbmc-machines/meta-openpower/common/recipes-phosphor/host/watchdog/openpower-host-watchdog-start@.target
[Unit]
Description=Start Host %I watchdog
Requires=phosphor-watchdog@%i.service
Conflicts=openpower-host-watchdog-stop@%i.target
RefuseManualStop=yes
