# meta-openbmc-machines/meta-openpower/common/recipes-phosphor/host/watchdog/openpower-host-watchdog-stop@.target
[Unit]
Description=Stop Host %I watchdog
Conflicts=phosphor-watchdog@%i.service
Conflicts=openpower-host-watchdog-start@%i.target
RefuseManualStop=yes
