# meta-phosphor/common/recipes-phosphor/dbus/obmc-mapper/obmc-mapper.service
[Unit]
Description=Phosphor OpenBMC DBus service management daemon

[Service]
Restart=always
Type=dbus
ExecStart=/usr/sbin/phosphor-mapper
BusName=org.openbmc.ObjectMapper

[Install]
WantedBy=multi-user.target
