# meta-quanta/meta-olympus-nuvoton/recipes-olympus-nuvoton/bmc-reboot-cause/files/bmc-reboot-cause.service
[Unit]
Description=BMC Reboot Cause service

[Service]
Restart=no
Type=oneshot
ExecStartPre=/usr/bin/env mapper wait /xyz/openbmc_project/object_mapper
ExecStartPre=/usr/bin/env mapper wait /xyz/openbmc_project/state/bmc0
ExecStart=/usr/bin/reboot-cause.sh
SyslogIdentifier=reboot-cause.sh

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
