# meta-facebook/meta-santabarbara/recipes-phosphor/mctp/files/mctp_setup.service
[Unit]
Description=MCTP Endpoint ID (EID) Setup
After=mctpd.service
Before=pldmd.service
Wants=mctpd.service
Requires=mctpd.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/mctp/mctp_setup.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target