# meta-facebook/meta-santabarbara/recipes-phosphor/mctp/files/mctp_setup@.service
[Unit]
Description=Setup %I MCTP Endpoint ID (EID)
After=mctpd.service
Before=pldmd.service
Before=xyz.openbmc_project.mctpreactor.service
Wants=mctpd.service
Requires=mctpd.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/mctp/mctp_setup.sh %I
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
