# meta/recipes-extended/at/at/atd.service
[Unit]
Description=Job spooling tools
After=syslog.target

[Service]
ExecStart=@SBINDIR@/atd -f

[Install]
WantedBy=multi-user.target
