# meta-qualcomm/meta-centriq2400-rep/recipes-phosphor/console/centriq-2400rep-console-client/centriq-2400rep-console-client@.service
[Unit]
Description=Control UART mux for obmc-console-client

[Service]
ExecStart=/usr/bin/env centriq-2400rep-console-client.sh
SyslogIdentifier=centriq-2400rep-console-client.sh

[Install]
WantedBy=multi-user.target
