# meta-phosphor/common/recipes-phosphor/mboxd/mboxd/mboxd.service
[Unit]
Description=Phosphor MBOX Daemon
Wants=dev-aspeed-mbox.device
After=dev-aspeed-mbox.device
Wants=dev-aspeed-lpc-ctrl.device
After=dev-aspeed-lpc-ctrl.device

[Service]
Restart=always
ExecStart={sbindir}/mboxd --flash {FLASH_SIZE} --verbose

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
