# meta-phosphor/common/recipes-phosphor/mboxd/mboxd/mboxd-reload@.service
[Unit]
Description=Reload mboxd during power on
Wants=obmc-power-start-pre@%i.target
Before=obmc-power-start-pre@%i.target
Conflicts=obmc-chassis-poweroff@%i.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart={sbindir}/mboxctl --reset

[Install]
WantedBy=obmc-chassis-poweron@%i.target
