# meta-phosphor/common/recipes-phosphor/mboxd/mboxd/mboxd-reload@.service
[Unit]
Description=Reload mboxd after power off
Wants=op-wait-power-off@%i.service
After=op-wait-power-off@%i.service
Conflicts=obmc-power-chassis-on@%i.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart={sbindir}/mboxctl --reset

[Install]
WantedBy=obmc-power-chassis-off@%i.target
WantedBy=obmc-quiesce-host@%i.target
