# meta-phosphor/common/recipes-phosphor/mboxd/mboxd/mboxd.service
[Unit]
Description=Phosphor MBOX Daemon

[Service]
Restart=always
ExecStart={sbindir}/mboxd
Wants=dev-aspeed-mbox.device
After=dev-aspeed-mbox.device

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
