# meta-phosphor/common/recipes-phosphor/mboxd/mboxd/org.openbmc.mboxd.service
[Unit]
Description=Phosphor IPMI MBOX Daemon

[Service]
Restart=always
ExecStart={sbindir}/mboxd
Type=dbus
BusName={BUSNAME}
Wants=dev-aspeed-mbox.device
After=dev-aspeed-mbox.device

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
