# meta-phosphor/common/recipes-phosphor/mboxd/mboxd/org.openbmc.mboxd.service
[Unit]
Description=Phosphor IPMI MBOX Daemon

[Service]
Restart=always
ExecStart={sbindir}/mboxd
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
