# meta-phosphor/common/recipes-phosphor/skeleton/skeleton/skeleton.service
[Unit]
Description="Temp placeholder for skeleton function"

[Service]
ExecStart=/usr/sbin/system_manager.py Palmetto

[Install]
WantedBy=multi-user.target
