# meta-openpower/recipes-phosphor/occ/openpower-occ-control/op-occ-disable@.service
[Unit]
Description=OpenPOWER OCC Active Disable
Before=obmc-host-stop-pre@%i.target
Before=op-enter-mpreboot@%i.service
Before=op-stop-instructions@%i.service
Conflicts=op-occ-enable@%i.service

[Service]
ExecStart=/usr/bin/occ-active.sh disable
Type=oneshot

[Install]
WantedBy=obmc-host-stop@%i.target
