# meta-openpower/recipes-phosphor/occ/openpower-occ-control/org.open_power.OCC.Control.service
[Unit]
Description=OpenPOWER OCC controller
Wants=mapper-wait@-xyz-openbmc_project-inventory.service
After=mapper-wait@-xyz-openbmc_project-inventory.service

[Service]
ExecStart=/usr/bin/env openpower-occ-control
SyslogIdentifier=openpower-occ-control
Restart=always
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy=multi-user.target
