# meta-facebook/meta-ventura2/recipes-phosphor/gpio/phosphor-gpio-monitor/gpio-set-high@.service
[Unit]
Description=%i GPIO Set High

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/gpio-setter "%i" 1
SyslogIdentifier=%i-gpio-set-high
