# meta-facebook/meta-ventura2/recipes-phosphor/gpio/phosphor-gpio-monitor/gpio-set-low@.service
[Unit]
Description=%i GPIO Set Low

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/gpio-setter "%i" 0
SyslogIdentifier=%i-gpio-set-low
