# meta-facebook/meta-ventura2/recipes-phosphor/gpio/phosphor-gpio-monitor/power-rail-assert-log@.service
[Unit]
Description=%i Power Rail Assert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/power-rail-event-logger "assert" "%i"
SyslogIdentifier=%i-power-rail-assert
