# meta-facebook/meta-sanmiguel/recipes-phosphor/gpio/phosphor-gpio-monitor/cpu-shdn-ok.service
[Unit]
Description=Cpu Shutdown OK

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/cpu-shdn-ok

SyslogIdentifier=cpu-shdn-ok
