# meta-facebook/meta-sanmiguel/recipes-phosphor/gpio/phosphor-gpio-monitor/cpu-boot-done-assert.service
[Unit]
Description=Cpu Boot Done Assert
Conflicts=cpu-boot-done-deassert.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/cpu-boot-done assert

SyslogIdentifier=cpu-boot-done-assert
