# meta-facebook/meta-sanmiguel/recipes-phosphor/gpio/phosphor-gpio-monitor/hmc-ready-assert.service
[Unit]
Description=HMC Ready Assert

[Service]
Type=simple
ExecStart=/usr/libexec/phosphor-gpio-monitor/hmc-ready assert

SyslogIdentifier=hmc-ready-assert
