# meta-facebook/meta-sanmiguel/recipes-phosphor/gpio/phosphor-gpio-monitor/prochot-assert-log@.service
[Unit]
Description=%i Prochot Assert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/thermal-event-logger "assert" "%i"
SyslogIdentifier=%i-prochot-assert
