# meta-facebook/meta-sanmiguel/recipes-phosphor/gpio/phosphor-gpio-monitor/rmc-main-power-assert.service
[Unit]
Description=RMC Main Power Disabled

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/rmc-main-power "disabled"

SyslogIdentifier=rmc-main-power-assert