# meta-facebook/meta-sanmiguel/recipes-phosphor/gpio/phosphor-gpio-monitor/rmc-main-power-deassert.service
[Unit]
Description=RMC Main Power Enabled

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/rmc-main-power "enabled"

SyslogIdentifier=rmc-main-power-deassert