# meta-facebook/meta-sanmiguel/recipes-phosphor/gpio/phosphor-gpio-monitor/thermtrip-assert-log@.service
[Unit]
Description=%i Thermtrip Assert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/thermal-event-logger "assert" "%i"
SyslogIdentifier=%i-thermtrip-assert
