# meta-quanta/meta-gsj/recipes-phosphor/usb-network/files/usb_network.service
[Unit]
Description=enable usb network

[Service]
Type=oneshot
ExecStart=/bin/bash -c '/usr/sbin/usb_network.sh'

[Install]
WantedBy=multi-user.target
