# meta-openbmc-machines/meta-openpower/meta-ibm/meta-witherspoon/recipes-phosphor/chassis/vrm-control/vrm-control@.service
[Unit]
Description=Apply voltage overrides to VRMs
Wants=avsbus-disable@%i.service
After=avsbus-disable@%i.service
Before=avsbus-enable@%i.service
Conflicts=obmc-chassis-poweroff@%i.target
ConditionPathExists=!/run/openbmc/chassis@%i-on

[Service]
ExecStart={bindir}/vrm-control.sh vdna=0.9 vdnb=0.9
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=obmc-chassis-poweron@%i.target
