# meta-facebook/meta-yosemite4/recipes-yosemite4/plat-svc/files/yosemite4-schematic-init.service
[Unit]
Description=Yosemite4 Schematic Init
After=xyz.openbmc_project.EntityManager.service
Requires=xyz.openbmc_project.EntityManager.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/yosemite4-schematic-init
SyslogIdentifier=Yosemite4 Schematic Init

[Install]
WantedBy=multi-user.target
