# meta-evb/meta-evb-nuvoton/meta-scm-npcm845/recipes-scm-npcm845/fan-init/fan-init/fan_init.service
[Unit]
Description=Fan init Daemon
Before=xyz.openbmc_project.EntityManager.service
Before=xyz.openbmc_project.fansensor.service

[Service]
Type=oneshot
ExecStart=/usr/bin/fan_init.sh

[Install]
WantedBy=multi-user.target
