# meta-phosphor/common/recipes-phosphor/inarp/inarp/inarp.service
[Unit]
Description=Phosphor Inverse ARP

[Service]
ExecStart={sbindir}/inarp
Restart=always

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
