# meta-phosphor/common/recipes-phosphor/inarp/inarp/inarp.service
[Unit]
Description=Phosphor Inverse ARP daemon

[Service]
ExecStart=/usr/sbin/inarp

[Install]
WantedBy=multi-user.target
