# meta-phosphor/common/recipes-phosphor/users/files/obmc-phosphor-user.service
[Unit]
Description=Phosphor OpenBMC user management daemon

[Service]
ExecStart={sbindir}/obmc-phosphor-userd
Environment="PYTHONUNBUFFERED=1"
Restart=always

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
