# common/recipes-phosphor/obmc-phosphor-sys/files/obmc-phosphor-sysd.service
[Unit]
Description=Phosphor OpenBMC system management daemon

[Service]
ExecStart=/usr/sbin/obmc-phosphor-sysd

[Install]
WantedBy=multi-user.target
