# meta-google/recipes-google/systemd/files/firmware-updates.target
[Unit]
Description=Firmware Updates and Configuration
Wants=firmware-updates-pre.target

[Install]
WantedBy=multi-user.target
