# meta-google/recipes-google/systemd/files/firmware-updates-pre.target
[Unit]
Description=Firmware Update Preconditions
