# meta-google/recipes-google/systemd/files/gbmc-psu-hardreset-pre.target
[Unit]
Description=Sys PSU Hard Reset
StopWhenUnneeded=true
Wants=gbmc-psu-hardreset-time.service
After=gbmc-psu-hardreset-time.service
Before=systemd-random-seed.service
