# meta-google/recipes-google/systemd/files/gbmc-log-ramoops.service
[Service]
Type=oneshot
SyslogIdentifier=gbmc-log-ramoops
SyslogLevelPrefix=true
ExecStart=/usr/bin/find /sys/fs/pstore -name '*-*' -print -exec cat {} \; -exec rm -f {} \;

[Install]
WantedBy=multi-user.target
