# meta-facebook/recipes-phosphor/console/obmc-console/fb-compute-multihost/obmc-console-ttyS1-ssh-host.socket
[Unit]
Description=Phosphor Host Console SSH socket for Slot2
Wants=obmc-console@ttyS1.service

[Socket]
ListenStream=2201
Accept=yes

[Install]
WantedBy=sockets.target
