# meta-facebook/recipes-phosphor/console/obmc-console/fb-compute-multihost/obmc-console-ttyS3-ssh-host.socket
[Unit]
Description=Phosphor Host Console SSH socket for Slot4
Wants=obmc-console@ttyS3.service

[Socket]
ListenStream=2203
Accept=yes

[Install]
WantedBy=sockets.target
