# meta-facebook/recipes-phosphor/console/obmc-console/fb-compute-multihost/obmc-console-ttyS1-ssh-bic.socket
[Unit]
Description=Phosphor BIC Console SSH socket for Slot2
Wants=obmc-console@ttyS1.service

[Socket]
ListenStream=2301
Accept=yes

[Install]
WantedBy=sockets.target
