# meta-facebook/recipes-phosphor/console/obmc-console/fb-compute-multihost/obmc-console-ttyS7-ssh-bic@.service
[Unit]
Description=Phosphor BIC Console SSH socket for Slot7
Wants=obmc-console@ttyS7.service

[Service]
ExecStartPre=/usr/bin/select-uart-mux slot7 bic
ExecStart=/usr/sbin/dropbear -i -r /etc/dropbear/dropbear_rsa_host_key -c "/usr/bin/obmc-console-client -i host6" $DROPBEAR_EXTRA_ARGS
SyslogIdentifier=console-client-ttyS7
ExecReload=/bin/kill -HUP $MAINPID
StandardInput=socket
KillMode=process
