# meta-facebook/recipes-phosphor/console/obmc-console/fb-compute-multihost/obmc-console-ttyS2-ssh-host.socket
[Unit]
Description=Phosphor Host Console SSH socket for Slot3
Wants=obmc-console@ttyS2.service

[Socket]
ListenStream=2202
Accept=yes

[Install]
WantedBy=sockets.target
