# meta-facebook/recipes-phosphor/console/obmc-console/fb-compute-multihost/obmc-console-ttyS0-ssh-bic.socket
[Unit]
Description=Phosphor BIC Console SSH socket for Slot1
Wants=obmc-console@ttyS0.service

[Socket]
ListenStream=2300
Accept=yes

[Install]
WantedBy=sockets.target
