# meta-openpower/common/recipes-phosphor/fsi/op-fsi/fsi-enable.service
[Unit]
Description=Enable FSI
Conflicts=fsi-disable.service
Wants=fsi-scan@0.service
Before=fsi-scan@0.service

[Service]
ExecStart=/bin/sh -c "echo 0 > /sys/devices/platform/gpio-fsi/external_mode"
Type=oneshot
RemainAfterExit=yes
