# meta-phosphor/recipes-phosphor/dbus/phosphor-mapper/mapper-subtree-remove@.service
[Unit]
Description=mapper subtree-remove %I
Wants=obmc-mapper.target
After=obmc-mapper.target

[Service]
Type=oneshot
RemainAfterExit=no
Restart=on-failure
TimeoutStartSec=infinity
ExecStart=/usr/bin/env mapper subtree-remove %I
SyslogIdentifier=mapper
