# meta-openpower/recipes-phosphor/debug/openpower-debug-collector/openpower-debug-collector-watchdog@.service
[Unit]
Description=OpenPOWER debug data collector for host watchdog timeout
Before=obmc-host-quiesce@%i.target
After=op-stop-instructions@%i.service

[Service]
ExecStart=/usr/bin/watchdog_timeout
Type=oneshot
Restart=no
