# yocto-poky/meta/recipes-extended/cups/cups/cups.service
[Unit]
Description=CUPS Printing Service

[Service]
ExecStart=@SBINDIR@/cupsd -f
PrivateTmp=true

[Install]
Also=cups.socket cups.path
WantedBy=printer.target
