# upstream-layers/openembedded-core/meta/recipes-kernel/lttng/lttng-tools/lttng-sessiond.service
[Unit]
Description=LTTng 2.x central tracing registry session daemon

[Service]
Type=forking
ExecStart=/usr/bin/lttng-sessiond -d

[Install]
WantedBy=multi-user.target
