# meta-phosphor/common/recipes-phosphor/obmc-phosphor-example-pydbus/files/obmc-phosphor-example-pydbus.service
[Unit]
Description=Phosphor OpenBMC QEMU application example

[Service]
Restart=always
ExecStart=/usr/sbin/obmc-phosphor-example-pydbus
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
