# meta-openembedded/meta-networking/recipes-protocols/rp-pppoe/rp-pppoe/pppoe-server.service
[Unit]
Description=PPPOE Service
After=network.target

[Service]
Type=forking
EnvironmentFile=-@SYSCONFDIR@/default/pppoe-server
ExecStart=@SBINDIR@/pppoe-server

[Install]
WantedBy=multi-user.target
