# meta-phosphor/common/recipes-phosphor/obmc-phosphor-sys/obmc-mgr-state/obmc-wait-standby.service
[Unit]
Description=Wait For Standby
Before=obmc-standby-reached.target

[Service]
Type=oneshot
ExecStart={base_bindir}/true
Restart=no

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
