# meta-fii/meta-kudo/recipes-phosphor/sensors/phosphor-virtual-sensor/cpu_detect_virtual.service
[Unit]
Description=Detect CPUs to change config
Requires=xyz.openbmc_project.adcsensor.service
After=xyz.openbmc_project.adcsensor.service
Before=phosphor-virtual-sensor.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-virtual-sensor/cpu_detect_virtual.sh

[Install]
WantedBy=multi-user.target
