# meta-phosphor/common/recipes-phosphor/dump/phosphor-debug-collector/var-lib-systemd-coredump.mount
[Unit]
Description=/var/lib/systemd/coredump

[Mount]
What=tmpfs
Where=/var/lib/systemd/coredump
Type=tmpfs
