# meta-phosphor/common/recipes-phosphor/dump/phosphor-debug-collector/xyz.openbmc_project.Dump.Monitor.service
[Unit]
Description=Phosphor Dump core monitor.

[Service]
ExecStart={sbindir}/phosphor-dump-monitor
Restart=always
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
