# meta-openembedded/meta-networking/recipes-daemons/vsftpd/files/vsftpd.service
[Unit]
Description=Vsftpd ftp daemon
After=network.target

[Service]
ExecStart=@SBINDIR@/vsftpd

[Install]
WantedBy=multi-user.target
