# meta-openembedded/meta-networking/recipes-support/ntp/ntp/ntpd.service
[Unit]
Description=Network Time Service
After=network.target

[Service]
Type=forking
PIDFile=/run/ntpd.pid
ExecStart=/usr/sbin/ntpd -u ntp:ntp -p /run/ntpd.pid -g

[Install]
WantedBy=multi-user.target
