# meta-quanta/meta-gbs/recipes-gbs/gbs-sysinit/files/gbs-sysinit.service
[Unit]
Description = Initialization for GBS boot up
Wants=mapper-wait@-xyz-openbmc_project-inventory.service
After=mapper-wait@-xyz-openbmc_project-inventory.service
Wants=mapper-wait@-xyz-openbmc_project-control-nvme.service
After=mapper-wait@-xyz-openbmc_project-control-nvme.service

[Service]
Type=oneshot
ExecStart=/usr/bin/gbs-sysinit.sh

[Install]
WantedBy=multi-user.target
