# meta-ampere/meta-jade/recipes-ampere/platform/mtjade-gpio-config/ampere-gpio-poweron@.service
[Unit]
Description=Ampere set GPIOs before power on
Before=op-wait-power-on@%i.service
After=obmc-power-start-pre@%i.target
Conflicts=obmc-chassis-poweroff@%i.target

[Service]
Type=oneshot
ExecStart=/usr/bin/env ampere_gpio_utils.sh power on
SyslogIdentifier=ampere_gpio_utils
