# meta-ampere/meta-jade/recipes-ampere/platform/mtjade-gpio-config/ampere-gpio-poweroff@.service
[Unit]
Description=Ampere set GPIOs before power off
Before=op-wait-power-off@%i.service
Conflicts=obmc-chassis-poweron@%i.target

[Service]
Type=oneshot
ExecStart=/usr/bin/env ampere_gpio_utils.sh power off
SyslogIdentifier=ampere_gpio_utils
