# meta-openbmc-machines/meta-openpower/meta-inventec/meta-lanyang/recipes-phosphor/chassis/set-vio/set-vio.service
[Unit]
Description=Set VIO rail voltages
Before={SYSTEMD_DEFAULT_TARGET}

[Service]
Type=oneshot

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
