# upstream-layers/openembedded-core/meta/recipes-core/busybox/files/busybox-klogd.service.in
[Unit]
Description=Kernel Logging Service

[Service]
ExecStart=@base_sbindir@/klogd -n

[Install]
WantedBy=multi-user.target
