# meta-facebook/meta-ventura/recipes-phosphor/leds/phosphor-led-manager/xyz.openbmc_project.LED.GroupManager.service.d/override.conf
[Unit]
Wants=phosphor-multi-gpio-monitor.service
Before=phosphor-multi-gpio-monitor.service
