# meta-phosphor/common/recipes-phosphor/state/discover-system-state/discover-system-state.target
[Unit]
Description=If power policy is ALWAYS_POWER_ON, powers the host on.
