# meta-phosphor/common/recipes-phosphor/state/discover-system-state/discover-system-state.target
# If power policy is ALWAYS_POWER_ON, powers the host on.
[Unit]
Description=Reboot If Enabled
