# meta-phosphor/recipes-phosphor/sensors/phosphor-hwmon/phosphor-max31785-msl@.service
[Unit]
Description=Phosphor MAX31785 MSL
ConditionFileNotEmpty={envfiledir}/obmc/hwmon-max31785/%I.conf
Wants=mapper-wait@-xyz-openbmc_project-inventory.service
After=mapper-wait@-xyz-openbmc_project-inventory.service

[Service]
Type=oneshot
ExecStart=/usr/bin/env max31785-msl -p "${{PATHS}}" -b ${{BUS}} -r ${{REVISION}}
EnvironmentFile={envfiledir}/obmc/hwmon-max31785/%I.conf
SyslogIdentifier=max31785-msl
RemainAfterExit=yes
