# meta-ibm/recipes-phosphor/bmcweb/bmcweb/bmcweb.socket
[Unit]
Description=Phosphor redfish socket (bmcweb)

[Socket]
ListenStream=127.0.0.1:8082

[Install]
WantedBy=sockets.target
