# upstream-layers/openembedded-core/meta/recipes-devtools/tcf-agent/tcf-agent/tcf-agent.service
[Unit]
Description=Target Communication Framework agent
After=network.target

[Service]
Type=forking
ExecStart=@SBINDIR@/tcf-agent -d -L- -l0
KillSignal=USR2
SuccessExitStatus=USR2

[Install]
WantedBy=multi-user.target
