# meta-phosphor/common/recipes-core/systemd/obmc-targets/obmc-host-reset-running@.target
[Unit]
Description=Host%i running after reset
Conflicts=obmc-host-shutdown@%i.target
RefuseManualStart=yes
RefuseManualStop=yes
