# meta-phosphor/common/recipes-core/systemd/obmc-targets/obmc-chassis-reset@.target
[Unit]
Description=Chassis%i (Reset Check)
Conflicts=obmc-power-chassis-off%i.target
RefuseManualStop=yes
