# meta-phosphor/common/recipes-core/systemd/obmc-targets/obmc-chassis@.target
[Unit]
Description=Chassis%i
After={SYSTEMD_DEFAULT_TARGET}
Requires={SYSTEMD_DEFAULT_TARGET}
