# meta-phosphor/common/recipes-core/systemd/obmc-targets/obmc-host-crash@.target
[Unit]
Description=Host instance %i crashed
Conflicts=obmc-chassis-poweroff@%i.target
Wants=obmc-host-quiesce@%i.target
RefuseManualStop=yes
