# meta-phosphor/common/recipes-core/systemd/obmc-targets/obmc-power-stop@.target
[Unit]
Description=Power%i Off
Wants=obmc-power-stop-pre@%i.target
After=obmc-power-stop-pre@%i.target
Wants={SYSTEMD_DEFAULT_TARGET}
Conflicts=obmc-chassis-poweron@%i.target
Conflicts=fsi-disable.service
RefuseManualStart=yes
RefuseManualStop=yes
