# meta-phosphor/common/recipes-core/systemd/obmc-targets/obmc-standby-reached.target
[Unit]
Description=Standby - Ready
RefuseManualStart=yes
RefuseManualStop=yes
