# meta-phosphor/common/recipes-core/systemd/obmc-targets/obmc-chassis-reset@.target
[Unit]
Description=Chassis%i (Reset Check)
Wants=mapper-wait@-org-openbmc-control-power%i.service
After=mapper-wait@-org-openbmc-control-power%i.service
Wants=mapper-wait@-xyz-openbmc_project-state-chassis%i.service
After=mapper-wait@-xyz-openbmc_project-state-chassis%i.service
Conflicts=obmc-chassis-stop@%i.target
RefuseManualStop=yes
OnFailureJobMode=flush

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
