# meta-phosphor/common/recipes-core/systemd/obmc-targets/obmc-fans-ready.target
[Unit]
Description=Fans Ready
RefuseManualStart=yes
RefuseManualStop=yes
OnFailure=obmc-quiesce-system@%i.target
OnFailureJobMode=flush