# meta-phosphor/common/recipes-core/systemd/obmc-targets/obmc-stop-host@.target
[Unit]
Description=Power%i Host Off
Requires={SYSTEMD_DEFAULT_TARGET}
Requires=op-stop-host@%i.service
After=op-stop-host@%i.service
RefuseManualStop=yes
OnFailure=obmc-quiesce-system@%i.target
OnFailureJobMode=flush