# meta-phosphor/common/recipes-core/systemd/obmc-targets/obmc-fan-control.target
[Unit]
Description=Fan Control
RefuseManualStart=yes
RefuseManualStop=yes
