# meta-phosphor/common/recipes-core/systemd/obmc-targets/obmc-host-timeout@.target
[Unit]
Description=Host instance %i timed out refreshing Watchdog
Conflicts=obmc-chassis-poweroff@%i.target
RefuseManualStop=yes
