# meta-phosphor/common/recipes-core/systemd/obmc-targets/obmc-power-reset-on@.target
[Unit]
Description=Chassis Power%i (On) after reset
Wants={SYSTEMD_DEFAULT_TARGET}
After={SYSTEMD_DEFAULT_TARGET}
Conflicts=obmc-power-chassis-off@%i.target
RefuseManualStart=yes
RefuseManualStop=yes