# meta-phosphor/common/recipes-core/systemd/obmc-targets/obmc-host-start-pre@.target
[Unit]
Description=Start Host%i (Pre)
Wants={SYSTEMD_DEFAULT_TARGET}
After={SYSTEMD_DEFAULT_TARGET}
Conflicts=obmc-host-stop@%i.target
RefuseManualStart=yes
RefuseManualStop=yes
