# meta-phosphor/common/recipes-core/systemd/obmc-targets/obmc-power-stop-pre@.target
[Unit]
Description=Power%i Off (Pre)
Wants={SYSTEMD_DEFAULT_TARGET}
After={SYSTEMD_DEFAULT_TARGET}
Conflicts=obmc-power-chassis-on@%i.target
RefuseManualStart=yes
RefuseManualStop=yes
