# meta-phosphor/common/recipes-core/systemd/obmc-targets/obmc-host-crash@.target
[Unit]
Description=Host instance %i crashed
Wants=obmc-host-started@%i.target
After=obmc-host-started@%i.target
Conflicts=obmc-chassis-poweroff@%i.target
RefuseManualStop=yes
