# meta-phosphor/common/recipes-core/systemd/obmc-targets/obmc-chassis-stop@.target
[Unit]
Description=Chassis%i (Stop)
Wants={SYSTEMD_DEFAULT_TARGET} mapper-wait@-org-openbmc-control-chassis%i.service
After={SYSTEMD_DEFAULT_TARGET} mapper-wait@-org-openbmc-control-chassis%i.service
Conflicts=obmc-chassis-start@%i.target
RefuseManualStop=yes
OnFailure=obmc-power-chassis-off@%i.target
OnFailureJobMode=flush
