# meta-phosphor/common/recipes-core/systemd/obmc-targets/obmc-mapper.target
[Unit]
Description=Phosphor Object Mapper
RefuseManualStart=yes
RefuseManualStop=yes
