# meta-phosphor/common/recipes-core/systemd/obmc-targets/obmc-gpio-asserted@.target
[Unit]
Description=GPIO asserted on %i
Wants=obmc-host-started@%i.target
After=obmc-host-started@%i.target
Conflicts=obmc-chassis-stop@%i.target
Conflicts=obmc-stop-host@%i.target
RefuseManualStart=yes
RefuseManualStop=yes
