# meta-phosphor/common/recipes-core/systemd/obmc-targets/obmc-host-stop-pre@.target
[Unit]
Description=Stop Host%i (Pre)
Wants={SYSTEMD_DEFAULT_TARGET}
After={SYSTEMD_DEFAULT_TARGET}
Conflicts=obmc-chassis-start@%i.target
RefuseManualStart=yes
RefuseManualStop=yes
