# meta-phosphor/common/recipes-core/systemd/obmc-targets/obmc-power-chassis-off@.target
[Unit]
Description=Chassis%i (Power Off)
Wants={SYSTEMD_DEFAULT_TARGET} mapper-wait@-xyz-openbmc_project-state-chassis%i.service
After={SYSTEMD_DEFAULT_TARGET} mapper-wait@-xyz-openbmc_project-state-chassis%i.service
Conflicts=obmc-power-chassis-on@%i.target
RefuseManualStop=yes
