# meta-phosphor/common/recipes-core/systemd/obmc-targets/obmc-standby.target
[Unit]
Description=Phosphor System - Standby
Requires=multi-user.target
AllowIsolate=yes

[Install]
WantedBy=multi-user.target
