# meta-phosphor/common/recipes-core/systemd/obmc-targets/obmc-fsi-disable.target
[Unit]
Description=Disable FSI Device Driver
