# meta-phosphor/common/recipes-core/systemd/obmc-targets/obmc-stop-host@.target
[Unit]
Description=Power%i Host Off
Wants={SYSTEMD_DEFAULT_TARGET}
RefuseManualStart=yes
RefuseManualStop=yes