# meta-phosphor/common/recipes-core/systemd/obmc-targets/obmc-power-reset-on@.target
[Unit]
Description=Chassis%i power on after reset
Conflicts=obmc-power-chassis-off@%i.target
RefuseManualStart=yes
RefuseManualStop=yes
