# meta-phosphor/common/recipes-core/systemd/obmc-targets/obmc-host-stopped@.target
[Unit]
Description=Host%i (Stopped)
Wants={SYSTEMD_DEFAULT_TARGET}
Wants=obmc-host-stopping@%i.target
After=obmc-host-stopping@%i.target
Before=obmc-power-stop-pre@%i.target
Conflicts=obmc-host-startmin@%i.target
RefuseManualStart=yes
RefuseManualStop=yes
