# meta-phosphor/common/recipes-core/systemd/obmc-targets/obmc-fan-watchdog-takeover.target
[Unit]
Description=Fan watchdog takeover target
