# meta-phosphor/common/recipes-core/systemd/obmc-targets/obmc-fan-control-ready@.target
[Unit]
Description=Fan Control Ready
RefuseManualStop=yes
Conflicts=obmc-chassis-poweroff@%i.target
