# meta-phosphor/common/recipes-core/systemd/obmc-targets/obmc-chassis-stop@.target
[Unit]
Description=Chassis%i (Stop)
Requires={SYSTEMD_DEFAULT_TARGET}
After={SYSTEMD_DEFAULT_TARGET}
Requires=mapper-wait@-org-openbmc-control-chassis%i.service
After=mapper-wait@-org-openbmc-control-chassis%i.service
Requires=obmc-power-chassis-off@%i.target obmc-stop-watchdog@%i.service
After=obmc-power-chassis-off@%i.target obmc-stop-watchdog@%i.service
Conflicts=obmc-chassis-start@%i.target
RefuseManualStop=yes
OnFailure=obmc-quiesce-host@%i.target
OnFailureJobMode=flush
