# meta-phosphor/common/recipes-core/systemd/obmc-targets/obmc-powered-on@.target
[Unit]
Description=Chassis%i Powered On
