# meta-phosphor/common/recipes-core/systemd/obmc-targets/reboot-guard.service
[Unit]
Description=Reboot Guard

[Service]
ExecStart=/bin/true

[Install]
RequiredBy=shutdown.target
