# meta-phosphor/common/recipes-core/systemd/obmc-targets/reboot-guard2.service
[Unit]
Description=Prevents BMC from rebooting.

[Service]
Type=oneshot
RemainAfterExit=true
ExecStop=<your script/program>

[Install]
WantedBy=multi-user.target
