# meta-phosphor/common/recipes-core/systemd/obmc-targets/start-reboot-guard.service
[Unit]
Description=Start Reboot Guard

[Service]
ExecStart=/bin/systemctl enable reboot-guard

[Install]
WantedBy=multi-user.target
