# meta-openembedded/meta-multimedia/recipes-multimedia/mycroft/files/mycroft-setup.service
[Unit]
Description=Mycroft Setup
DefaultDependencies=no
After=systemd-user-sessions.service

[Service]
Type=oneshot
ExecStart=cp -rf @LIBDIR@/mycroft /var/
ExecStart=/var/mycroft/dev_setup.sh --allow-root

[Install]
WantedBy=multi-user.target
