# meta-openembedded/meta-multimedia/recipes-multimedia/mycroft/files/mycroft.service
[Unit]
Description=Mycroft
RequiresMountsFor=/var
After=mycroft-setup.service

[Service]
Type=forking
WorkingDirectory=@LIBDIR@/mycroft
ExecStartPre=mkdir -p /var/log/mycroft
ExecStart=@LIBDIR@/mycroft/start-mycroft.sh all
ExecStop=@LIBDIR@/mycroft/stop-mycroft.sh all

[Install]
WantedBy=multi-user.target
