# meta-openembedded/meta-multimedia/recipes-multimedia/mycroft/files/mycroft.service
[Unit]
Description=Mycroft
RequiresMountsFor=/var
After=mycroft-setup.service

[Service]
Type=forking
WorkingDirectory=/var/mycroft
ExecStartPre=mkdir -p /var/log/mycroft
ExecStart=/var/mycroft/start-mycroft.sh all
ExecStop=/var/mycroft/stop-mycroft.sh all

[Install]
WantedBy=multi-user.target
