# meta-quanta/meta-olympus-nuvoton/recipes-phosphor/smbios/smbios-mdr/xyz.openbmc_project.cpuinfo.service
[Unit]
Description=Intel BMC CPU Info
Requires=smbios-mdrv2.service
After=smbios-mdrv2.service

[Service]
Restart=always
ExecStart=/usr/bin/cpuinfoapp
Type=dbus
BusName=xyz.openbmc_project.CPUInfo
Nice=19

[Install]
WantedBy=multi-user.target
