# meta-phosphor/recipes-phosphor/ipmi/phosphor-ipmi-bt/org.openbmc.HostIpmi.service
[Unit]
Description=Phosphor IPMI BT DBus Bridge

[Service]
Restart=always
ExecStart=/usr/bin/env btbridged
SyslogIdentifier=btbridged
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy=multi-user.target
