# meta-openbmc-machines/meta-arm/meta-qualcomm/meta-centriq2400-rep/recipes-phosphor/chassis/ir38163-workaround/ir-add@.service
[Unit]
Description=Add IR38163 device
Wants=obmc-power-start@%i.target
Before=obmc-power-start@%i.target
Conflicts=obmc-chassis-poweroff@%i.target
ConditionPathExists=!/run/openbmc/chassis@%i-on

[Service]
ExecStart=/usr/bin/env ir38163_workaround.sh add
SyslogIdentifier=ir38163_workaround.sh
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=obmc-host-start@%i.target

