# meta-phosphor/common/recipes-phosphor/host/phosphor-host-postd/lpcsnoop.service
[Unit]
Description=LPC Snoop Daemon
After=systemd-modules-load.service

[Service]
Restart=always
ExecStart={sbindir}/snoopd

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
