# upstream-layers/meta-openembedded/meta-oe/recipes-support/openldap/openldap/slapd.service
[Unit]
Description=Standalone LDAP Daemon
After=syslog.target network.target

[Service]
Type=forking
ExecStart=@SBINDIR@/slapd

[Install]
WantedBy=multi-user.target
