# meta-phosphor/common/recipes-phosphor/obmc-phosphor-sys/obmc-control-bmc/org.openbmc.control.Bmc@.service
[Unit]
Description=Phosphor BMC%i Control

[Service]
Restart=always
ExecStart={sbindir}/control_bmc.exe
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
