# meta-ampere/meta-common/recipes-phosphor/gpio/phosphor-gpio-monitor/ampere_hightemp_start@.service
[Unit]
Description=Ampere Computing LLC HighTempt event

[Service]
Restart=no
ExecStart=/usr/sbin/ampere_add_redfishevent.sh OpenBMC.0.1.AmpereWarning "CPU%i HighTemp asserted"
