# meta-ampere/meta-common/recipes-phosphor/gpio/phosphor-gpio-monitor/ampere_psu_power_restore@.service
[Unit]
Description=Ampere Computing LLC PSU Power Restored event
StartLimitIntervalSec=0

[Service]
Restart=no
ExecStart=/usr/sbin/ampere_add_redfishevent.sh OpenBMC.0.1.PowerSupplyPowerRestored "PSU%i"
