# meta-ampere/meta-common/recipes-phosphor/gpio/phosphor-gpio-monitor/ampere_hightemp_start@.service
[Unit]
Description=Ampere Computing LLC HighTemp Assertion event
StartLimitIntervalSec=0

[Service]
Restart=no
ExecStart=/usr/sbin/ampere_hightemp.sh %i asserted
