# meta-ampere/meta-common/recipes-phosphor/gpio/phosphor-gpio-monitor/ampere_overtemp@.service
[Unit]
Description=Ampere Computing LLC Overtemp event handling

[Service]
Restart=no
ExecStart=/usr/sbin/ampere_overtemp.sh %i
