# meta-ampere/meta-common/recipes-phosphor/gpio/phosphor-gpio-monitor/ampere_hightemp_stop@.service
[Unit]
Description=Ampere Computing LLC HighTemp Deassertion event
StartLimitIntervalSec=0

[Service]
Restart=no
ExecStart=/usr/sbin/ampere_add_redfishevent.sh OpenBMC.0.1.AmpereEvent "CPU%i HighTemp deasserted"
