# meta-phosphor/common/recipes-phosphor/users/obmc-phosphor-user/org.openbmc.UserManager.service
[Unit]
Description=Phosphor User Manager

[Service]
ExecStart=/usr/bin/env obmc-phosphor-userd
SyslogIdentifier=obmc-phosphor-userd
Environment="PYTHONUNBUFFERED=1"
Restart=always
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
