# meta-openpower/recipes-phosphor/host/hostboot-settings/hostboot-settings.service
[Unit]
Description=Hostboot Boot Settings

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/bin/hb_settings --uart vuart

[Install]
WantedBy=multi-user.target
