# import-layers/meta-openembedded/meta-gnome/recipes-gnome/gdm/gdm/gdm.service.in
[Unit]
Description=Gnome Display Manager
After=systemd-user-sessions.service

[Service]
ExecStart=%sbindir%/gdm -nodaemon
StandardOutput=syslog

[Install]
Alias=display-manager.service
