# meta-fii/meta-mori/recipes-phosphor/ipmi/phosphor-ipmi-flash/bios-update/phosphor-ipmi-flash-bios-update.service
[Unit]
Description=Phosphor-ipmi-flash update BIOS firmware service

[Service]
ExecStart=/usr/bin/mori-fw.sh bios /run/initramfs/bios-image
Type=oneshot
