# meta-openembedded/meta-oe/recipes-support/libusbgx/libusbgx/usbgx.service
[Unit]
Description=Load USB gadget schemas

[Service]
Type=oneshot
ExecStart=/usr/bin/gadget-start

[Install]
WantedBy=multi-user.target
