# meta-phosphor/common/recipes-phosphor/ipmi/phosphor-ipmi-net/phosphor-ipmi-net.service
[Unit]
Description=Network IPMI daemon
Wants=mapper-wait@-.service
After=mapper-wait@-.service

[Service]
ExecStart={sbindir}/netipmid
Restart=always
