# meta-phosphor/common/recipes-phosphor/ipmi/phosphor-ipmi-net/phosphor-ipmi-net.service
[Unit]
Description=Network IPMI daemon

[Service]
ExecStart={sbindir}/netipmid
Restart=always
