# meta-facebook/meta-bletchley/recipes-bletchley/motor-ctrl/files/host-poweron@.service
[Unit]
Description=Power on System in Sled%i by Step Motor
Requires=motor-init-calibration@%i.service
After=motor-init-calibration@%i.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/power-ctrl sled%i on
SyslogIdentifier=power-ctrl

