# meta-facebook/meta-bletchley/recipes-bletchley/motor-ctrl/files/host-ac-on@.service
[Unit]
Description=Turn on AC at Sled%i
After=motor-init-calibration@%i.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/power-ctrl sled%i ac-on
SyslogIdentifier=power-ctrl

[Install]
WantedBy=multi-user.target
