# meta-facebook/meta-bletchley/recipes-bletchley/motor-ctrl/files/host-ac-off@.service
[Unit]
Description=Turn off AC at Sled%i
Conflicts=host-ac-on@%i.service
After=host-poweroff@%i.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/power-ctrl sled%i ac-off
SyslogIdentifier=power-ctrl
