# meta-openbmc-machines/meta-openpower/meta-ingrasys/meta-zaius/recipes-phosphor/chassis/refclock-fixup/op-refclock-fixup@.service
[Unit]
Description=Fix up Zaius host refclock
Wants=obmc-host-start-pre@%i.target
Before=obmc-host-start-pre@%i.target
Wants=obmc-power-on@%i.target
After=obmc-power-on@%i.target
Before=fsi-scan@%i.service
Conflicts=obmc-host-stop@%i.target
ConditionPathExists=!/run/openbmc/chassis@%i-on

[Service]
Restart=no
Type=oneshot
RemainAfterExit=yes
ExecStart={sbindir}/fix_zaius_refclock.sh

[Install]
RequiredBy=obmc-chassis-poweron@%i.target
