# meta-ibm/recipes-httpd/nginx/files/nginx.socket
[Unit]
Description=Nginx

[Socket]
ListenStream=8081

[Install]
WantedBy=sockets.target
