# meta-phosphor/recipes-phosphor/system/obmc-phosphor-sysd/obmc-phosphor-sysd.service
[Unit]
Description=Phosphor System Manager

[Service]
Restart=always
ExecStart=/usr/bin/env obmc-phosphor-sysd
SyslogIdentifier=obmc-phosphor-sysd

[Install]
WantedBy=multi-user.target
