# meta-ufispace/meta-ncplite/recipes-phosphor/state/phosphor-state-manager/ncplite-chassis-poweron@.service
[Unit]
Description=NCPLite Chassis Power on
Before=phosphor-wait-power-on@%i.service
After=obmc-power-start-pre@%i.target
Conflicts=ncplite-chassis-poweroff@%i.target

[Service]
Type=oneshot
ExecStart=@SBINDIR@/chassis-control on
SyslogIdentifier=chassis-control
