# poky/meta/recipes-kernel/kexec/kexec-tools/kdump.service
[Unit]
Description=Reboot and dump vmcore via kexec
DefaultDependencies=no
Requires=sysinit.target
After=sysinit.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@LIBEXECDIR@/kdump-helper start
ExecStop=@LIBEXECDIR@/kdump-helper stop

[Install]
WantedBy=multi-user.target
