# meta-openembedded/meta-networking/recipes-filter/ebtables/ebtables-2.0.10-4/ebtables.service
[Unit]
Description=Ethernet Bridge Filtering Tables

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@SBINDIR@/ebtables.common start
ExecStop=@SBINDIR@/ebtables.common stop

[Install]
WantedBy=multi-user.target
