# import-layers/meta-openembedded/meta-oe/recipes-extended/openwsman/openwsman/openwsmand.service
[Unit]
Description=Openwsman WS-Management Service
After=syslog.target

[Service]
Type=forking
ExecStart=@SBINDIR@/openwsmand
ExecStartPre=@SYSCONFDIR@/openwsman/owsmangencert.sh
PIDFile=@LOCALSTATEDIR@/run/wsmand.pid

[Install]
WantedBy=multi-user.target
