# meta-fii/meta-mori/recipes-mori/hotswap-power-cycle/files/hotswap-power-cycle.service
[Unit]
Description=Power Cycle by Hotswap Controller
After=gbmc-psu-hardreset.target

[Service]
Type=oneshot
ExecStart=gpioset 2 30=1
StandardOutput=syslog

[Install]
WantedBy=gbmc-psu-hardreset.target
