# meta-phosphor/common/recipes-phosphor/fans/files/obmc-phosphor-fand.service
[Unit]
Description=Phosphor OpenBMC fan management daemon

[Service]
Restart=always
ExecStart={sbindir}/obmc-phosphor-fand

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
