# meta-phosphor/common/recipes-phosphor/fans/files/obmc-phosphor-fand.service
[Unit]
Description=Phosphor Fan Controller

[Service]
Restart=always
ExecStart=/usr/bin/env obmc-phosphor-fand
SyslogIdentifier=obmc-phosphor-fand

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
