# meta-phosphor/common/recipes-phosphor/inventory/phosphor-inventory-manager/xyz.openbmc_project.Inventory.Manager.service
[Unit]
Description=Phosphor Inventory Manager

[Service]
Restart=always
Type=dbus
ExecStart={sbindir}/phosphor-inventory
BusName={BUSNAME}

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
