# meta-openbmc-machines/meta-openpower/meta-ibm/meta-romulus/recipes-phosphor/chassis/vddr-control/vddr-control@.service
[Unit]
Description=Apply voltage overrides to VDDRs
Wants=cpld_trigger@%i.service
After=cpld_trigger@%i.service
Wants=obmc-host-start-pre@%i.target
Before=obmc-host-start-pre@%i.target
Conflicts=obmc-power-chassis-off@%i.target

[Service]
ExecStart={bindir}/vddr.sh
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=obmc-power-chassis-on@%i.target
