# poky/meta/recipes-support/rng-tools/rng-tools/rngd.service
[Unit]
Description=Hardware RNG Entropy Gatherer Daemon
DefaultDependencies=no
After=systemd-udev-settle.service
Before=sysinit.target

[Service]
ExecStart=@SBINDIR@/rngd -f -r /dev/hwrng
SuccessExitStatus=66

[Install]
WantedBy=multi-user.target
