# meta-openembedded/meta-networking/recipes-daemons/cyrus-sasl/cyrus-sasl/saslauthd.service
[Unit]
Description=SASL authentication daemon.
After=syslog.target

[Service]
Type=forking
PIDFile=/run/saslauthd/saslauthd.pid
EnvironmentFile=@SYSCONFDIR@/default/saslauthd
ExecStart=@SBINDIR@/saslauthd -m $SOCKETDIR -a $MECH $FLAGS

[Install]
WantedBy=multi-user.target
