# meta-phosphor/common/recipes-phosphor/obmc-phosphor-host/obmc-op-control-host/org.openbmc.control.Host@.service
[Unit]
Description=OpenPOWER Host%i Control
Wants=mapper-wait@-org-openbmc-managers-System.service
After=mapper-wait@-org-openbmc-managers-System.service

[Service]
Restart=always
ExecStart={sbindir}/control_host.exe
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
