# meta-openbmc-machines/meta-openpower/common/recipes-phosphor/host/checkstop-monitor/openpower-host-checkstop-monitor-start@.target
[Unit]
Description=Start Host %I monitor
Requires=phosphor-gpio-monitor@%i.service
RefuseManualStop=yes
