# meta-openbmc-machines/meta-openpower/common/recipes-phosphor/host/checkstop-monitor/openpower-host-checkstop-monitor-stop@.target
[Unit]
Description=Stop Host %I monitor
Conflicts=phosphor-gpio-monitor@%i.service
RefuseManualStop=yes
