# meta-evb/meta-evb-nuvoton/meta-scm-npcm845/recipes-scm-npcm845/usb-loopback/files/usb_loopback.service
[Unit]
Description= USB loopback test

[Service]
Type=simple
RemainAfterExit=no
ExecStart=/bin/bash -c '/usr/sbin/usb_loopback.sh /dev/sda /tmp/mnt'