# meta-evb/meta-evb-nuvoton/meta-buv-runbmc/recipes-buv-runbmc/power/first-boot-set-psu/first-boot-set-psu@.service
[Unit]
Description=Init PSU Inventory for PSU %I
Wants=mapper-wait@-xyz-openbmc_project-inventory.service
After=mapper-wait@-xyz-openbmc_project-inventory.service

[Service]
ExecStart=/usr/bin/first-boot-set-psu.sh %i
RemainAfterExit=yes
SyslogIdentifier=first-boot-set-psu
Restart=on-failure
RestartSec=20

[Install]
WantedBy=multi-user.target
