# meta-quanta/meta-olympus-nuvoton/recipes-phosphor/chassis/olympus-nuvoton-powerctrl/host-gpio.service
[Unit]
Description=RunBMC Configure GPIOs for Host Power Control

[Service]
Restart=no
RemainAfterExit=true
Type=oneshot
ExecStart=/usr/bin/env init_once.sh
SyslogIdentifier=init_once.sh

[Install]
WantedBy=multi-user.target
