# common/recipes-phosphor/obmc-phosphor-chassis/files/obmc-phosphor-chassisd.service
[Unit]
Description=Phosphor OpenBMC chassis management daemon

[Service]
ExecStart=/usr/sbin/obmc-phosphor-chassisd

[Install]
WantedBy=multi-user.target
