# import-layers/meta-openembedded/meta-oe/recipes-support/openipmi/files/ipmi.service
[Unit]
Description=IPMI Driver
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@LIBEXECDIR@/openipmi-helper start
ExecStop=@LIBEXECDIR@/openipmi-helper stop-all

[Install]
WantedBy=multi-user.target
