# meta-hpe/meta-rl300-g11/recipes-hpe/power-sequencing/files/gpiosdown.service
[Unit]
Description=HPE RL300 PowerOK

[Service]
ExecStart=/usr/bin/gpios-manager down

[Install]
# WantedBy=multi-user.target
