# meta-hpe/meta-rl300-g11/recipes-hpe/power-sequencing/files/powerbuttonreleased.service
[Unit]
Description=HPE RL300 Power Button Release

[Service]
ExecStart=/usr/bin/power-button released

[Install]
# WantedBy=multi-user.target
