# meta-hpe/meta-rl300-g11/recipes-hpe/power-sequencing/files/powerbuttonpressed.service
[Unit]
Description=HPE RL300 Power Button

[Service]
ExecStart=/usr/bin/power-button pressed

[Install]
# WantedBy=multi-user.target
