# meta-intel/meta-common/recipes-intel/dbus-sensors/dbus-sensors/xyz.openbmc_project.cpusensor.service
[Unit]
Description=CPU Sensor
StopWhenUnneeded=false

[Service]
Restart=always
RestartSec=5
StartLimitBurst=10
ExecStart={bindir}/cpusensor

[Install]
WantedBy=multi-user.target
