# meta-intel/meta-common/recipes-intel/dbus-sensors/dbus-sensors/xyz.openbmc_project.adcsensor.service
[Unit]
Description=Adc Sensor
StopWhenUnneeded=false

[Service]
Restart=always
RestartSec=5
ExecStart={bindir}/adcsensor

[Install]
WantedBy=multi-user.target
