# meta-openbmc-machines/meta-openpower/meta-ibm/recipes-phosphor/phosphor-gevent/phosphor-gevent/phosphor-gevent.service
[Unit]
Description=Phosphor REST Server
After=network.target
After=obmc-webserver-pre.target

[Service]
Restart=always
ExecStart=/usr/bin/env phosphor-gevent $APPLICATION --no-ssl
SyslogIdentifier=phosphor-gevent
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile={envfiledir}/obmc/wsgi_app
