# meta-openbmc-machines/meta-openpower/meta-ibm/recipes-phosphor/phosphor-gevent/phosphor-gevent/phosphor-gevent.socket
[Unit]
Description=Phosphor Webserver socket

[Socket]
ListenStream=127.0.0.1:8081

[Install]
WantedBy=sockets.target
