# import-layers/meta-openembedded/meta-networking/recipes-daemons/tftp-hpa/files/tftpd-hpa.service
[Unit]
Description=Tftp Server
Requires=tftpd-hpa.socket

[Service]
ExecStart=@SBINDIR@/in.tftpd-hpa -s /var/lib/tftpboot
StandardInput=socket

[Install]
Also=tftpd-hpa.socket
