# meta-facebook/meta-yosemite4/recipes-phosphor/pldm/files/nic-powercycle@.service
[Unit]
Description=Nic power cycle

[Service]
Type=oneshot
ExecStart=/usr/libexec/pldm/nic-powercycle %i
