# meta-facebook/meta-yosemite4/recipes-phosphor/pldm/files/pldmd-exit-handler.service
[Unit]
Description=PLDMD Exit Handler

[Service]
Type=oneshot
ExecStart=/usr/libexec/pldm/pldmd-exit-handler
SyslogIdentifier=pldmd-exit-handler
