# meta-facebook/meta-yosemite4/recipes-phosphor/pldm/files/pldm_nic_power_cycle@.service
[Unit]
Description=PLDM nic power cycle

[Service]
Type=oneshot
ExecStart=/usr/libexec/pldm/pldm_nic_power_cycle %i
