# import-layers/meta-openembedded/meta-oe/recipes-connectivity/zabbix/zabbix/zabbix-agent.service
[Unit]
Description=Zabbix Monitor Agent
After=syslog.target network.target

[Service]
Type=oneshot
ExecStart=@SBINDIR@/zabbix_agentd
RemainAfterExit=yes
User=zabbix

[Install]
WantedBy=multi-user.target
