# meta-google/recipes-phosphor/flash/dummy-gbmc-update/dummy-verify.service
[Unit]
Description=Dummy flash file verification

[Service]
Type=oneshot
ExecStart=/bin/mv /run/initramfs/bmc-image /run/initramfs/dummy
