# meta-phosphor/common/recipes-phosphor/examples/obmc-phosphor-example-sdbus/org.openbmc.examples.SDBusService.service
[Unit]
Description=Phosphor Application Example

[Service]
Restart=always
ExecStart={sbindir}/obmc-phosphor-example-sdbus

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
