# meta-hpe/meta-dl385-g11/recipes-hpe/power-supplies-monitoring/files/psu1up.service
[Unit]
Description=HPE DL360 PSU1 Manager

[Service]
ExecStart=/usr/bin/psus-manager.sh 1 up

[Install]
# WantedBy=multi-user.target
