# meta-ampere/meta-jade/recipes-ampere/host/ampere-scp-failover/ampere-scp-failover.service
[Unit]
Description=Ampere Computing LLC SCP failover service

[Service]
Restart=always
ExecStart=/usr/bin/ampere-scp-failover
SyslogIdentifier=ampere-scp-failover

[Install]
WantedBy=multi-user.target
