# meta-phosphor/recipes-phosphor/video/uart-render-controller/uart-render-controller.service

[Unit]
Description=UART render controller
Requires=xyz.openbmc_project.State.Host.service
After=xyz.openbmc_project.State.Host.service

[Service]
ExecStart=/usr/sbin/obmc-uart-render-controller \
	  /sys/devices/platform/ahb/ahb:apb/1e6e2000.syscon/1e6e2000.syscon:scratch/vga0 \
	  /sys/devices/platform/ahb/ahb:apb/1e6e2000.syscon/1e6e2000.syscon:scratch/dac_mux \
	  fbterm.service

[Install]
WantedBy=multi-user.target
