# meta-ufispace/meta-ncplite/recipes-ncplite/ncplite-inventory-log/files/inventory-log.service
[Unit]
Description=Check GPIO State Manager
After=phosphor-gpio-presence@.service
Wants=phosphor-gpio-presence@.service

[Service]
ExecStart=@LIBEXECDIR@/ncplite-inventory-log/inventory-log.sh
SyslogIdentifier=inventory-log
Type=simple

[Install]
WantedBy=multi-user.target
