# meta-openembedded/meta-oe/recipes-devtools/android-tools/android-tools/android-tools-adbd.service
[Unit]
Description=Android Debug Bridge

[Service]
Type=simple
Restart=on-failure
ExecStartPre=/usr/bin/android-gadget-setup adb
ExecStart=/usr/bin/adbd
StandardOutput=null

[Install]
WantedBy=basic.target
