# meta-phosphor/recipes-phosphor/host/obmc-op-control-host/org.openbmc.control.Host@.service
[Unit]
Description=OpenPOWER Host%i Control

[Service]
Restart=always
ExecStart=/usr/bin/env control_host.exe
SyslogIdentifier=control_host.exe
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
