# meta-openbmc-machines/meta-openpower/meta-ibm/meta-witherspoon/recipes-phosphor/chassis/vcs-pre-post/vcs_off@.service
[Unit]
Description=Disable VCS before Witherspoon power on
Wants=mapper-wait@-org-openbmc-control-power%i.service
Before=mapper-wait@-org-openbmc-control-power%i.service

[Service]
Restart=always
ExecStart={bindir}/vcs_off.sh
Type=simple
RemainAfterExit=yes

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
