# meta-openbmc-machines/meta-openpower/meta-ibm/meta-witherspoon/recipes-phosphor/chassis/vcs-pre-post/vcs_on@.service
[Unit]
Description=Turn VCS on
Wants=mapper-wait@-vcs_workaround%i.service
After=mapper-wait@-vcs_workaround%i.service

[Service]
Restart=always
ExecStart={bindir}/vcs_on.sh
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}

