# meta-openbmc-machines/meta-openpower/meta-ibm/meta-witherspoon/recipes-phosphor/chassis/vcs-pre-post/vcs_on@.service
[Unit]
Description=Turn VCS on
Wants=vcs_workaround@.service
After=vcs_workaround@.service
Wants=obmc-host-start@%i.target
Before=obmc-host-start@%i.target
Conflicts=obmc-chassis-stop@%i.target

[Service]
ExecStart={bindir}/vcs_on.sh
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=obmc-chassis-start@%i.target
