# meta-phosphor/recipes-phosphor/flash/obmc-flash-bios/org.openbmc.control.Flash.service
[Unit]
Description=Phosphor BIOS Code Update
Wants=mapper-wait@-org-openbmc-managers-Download.service
After=mapper-wait@-org-openbmc-managers-Download.service

[Service]
Restart=always
ExecStart=/usr/bin/env flash_bios.exe
SyslogIdentifier=flash_bios.exe
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy=multi-user.target
