# meta-phosphor/common/recipes-phosphor/events/files/obmc-phosphor-event.service
[Unit]
Description=Phosphor Event Manager

[Service]
EnvironmentFile=-/etc/default/eventd/eventd.conf
Restart=always
ExecStart={sbindir}/obmc-phosphor-eventd $OPTIONS

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
