# meta-quanta/meta-gbs/recipes-x86/chassis/x86-power-control/chassis-system-reset.service
[Unit]
Description=System unit to hard reset or system reset chassis

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/tray_powercycle.sh

[Install]
WantedBy=chassis-system-reset.target
