# meta-phosphor/common/recipes-phosphor/flash/files/obmc-phosphor-flashd.service
[Unit]
Description=Phosphor OpenBMC flash management daemon

[Service]
Restart=always
ExecStart=/usr/sbin/obmc-phosphor-flashd

[Install]
WantedBy=multi-user.target
