# meta-hpe/meta-dl360-g11/recipes-hpe/power-sequencing/files/gpiosup.service
[Unit]
Description=HPE DL360 PowerOK

[Service]
ExecStart=/usr/bin/gpios-manager.sh up

[Install]
# WantedBy=multi-user.target
