# meta-hpe/meta-dl360-g11/recipes-hpe/power-sequencing/files/gpios-manager.service
[Unit]
Description=HPE virtualnic configuration
Before=network-pre.target
Wants=network-pre.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/etc/virtualnic/virtualnic start
ExecStop=/etc/virtualnic/virtualnic stop

[Install]
WantedBy=multi-user.target
