# meta-hpe/meta-dl360-g11/recipes-hpe/power-sequencing/files/monitor-ipmi-power-transition.service
[Unit]
Description=DL360 IPMI power sequencing request
After=dbus.service

[Service]
Restart=always
RestartSec=5
ExecStart=/usr/bin/ipmi-power-button-monitoring.sh

[Install]
WantedBy=multi-user.target
