# meta-hpe/meta-dl360-g11/recipes-hpe/power-sequencing/files/powerbuttonreleased.service
[Unit]
Description=HPE DL360 Power Button Release

[Service]
ExecStart=/usr/bin/power-button.sh released

[Install]
# WantedBy=multi-user.target
