# meta-nvidia/recipes-nvidia/usb-ethernet-init/files/usb-ethernet-init.service
[Unit]
Description=Init USB Host Ethernet Gadget
Before=network.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/usb-ethernet-init.sh

[Install]
WantedBy=multi-user.target
