# import-layers/meta-openembedded/meta-networking/recipes-support/geoip/geoipupdate-2.5.0/geoipupdate.cron
# top of crontab
MAILTO=your@email.com

32 11 * * 4 /usr/local/bin/geoipupdate
# end of crontab
