# meta-fii/meta-mori/recipes-phosphor/dbus/led-policy-mori/boot_status_led_on.service
[Unit]
Description=Mori boot status led on service

[Service]
Type=oneshot
ExecStart=/bin/sh -c "busctl set-property xyz.openbmc_project.LED.GroupManager /xyz/openbmc_project/led/groups/boot_status_inactive xyz.openbmc_project.Led.Group Asserted b false && \
busctl set-property xyz.openbmc_project.LED.GroupManager /xyz/openbmc_project/led/groups/boot_status_bootcomplete xyz.openbmc_project.Led.Group Asserted b true"
SyslogIdentifier=phosphor-dbus-monitor
