# meta-evb/meta-evb-nuvoton/meta-scm-npcm845/recipes-scm-npcm845/dhcp-check/dhcp-check/dhcp-check.service
[Unit]
Description=Check for DHCP address
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/dhcp-check.sh

[Install]
WantedBy=multi-user.target