# meta-ampere/meta-mitchell/recipes-ampere/platform/ampere-utils/ampere-fanctrl-manual.service
[Unit]
Description= Ampere Manual Fan Control Service

[Service]
ExecStart=echo "Disabling Ampere Fan Control"
Type=oneshot
RemainAfterExit=yes