# meta-phosphor/recipes-phosphor/flash/phosphor-software-manager/obmc-flash-bmc-ubiro-remove@.service
[Unit]
Description=Deletes read-only and kernel ubi volume %I

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/env obmc-flash-bmc ubiremove rofs-%i
ExecStart=/usr/bin/env obmc-flash-bmc ubiremove kernel-%i
SyslogIdentifier=obmc-flash-bmc
