# meta-phosphor/recipes-phosphor/flash/phosphor-software-manager/xyz.openbmc_project.Software.BMC.Updater.service
[Unit]
Description=OpenBMC Software Update Manager
Wants=xyz.openbmc_project.Software.Version.service
Before=xyz.openbmc_project.Software.Version.service
Wants=obmc-mapper.target
After=obmc-mapper.target

[Service]
ExecStart=/usr/bin/env phosphor-image-updater
SyslogIdentifier=phosphor-image-updater
Restart=always
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
