# meta-facebook/recipes-phosphor/state/phosphor-state-manager/greatlakes/host-poweron@.service
[Unit]
Description=power on host:%i

[Service]
Type=oneshot
ExecStart=/usr/libexec/host-poweron %i

[Install]
WantedBy=obmc-host-start@%i.target
