# meta-facebook/recipes-phosphor/state/phosphor-state-manager/greatlakes/power-ctrl-init.service
[Unit]
Description=Power Control Initialization
Before=mapper-wait@-xyz.openbmc_project.State.Chassis@0.service
Before=mapper-wait@-xyz-openbmc_project-state-bmc.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-state-manager/power-ctrl-init

[Install]
WantedBy=multi-user.target
