# meta-facebook/recipes-phosphor/state/phosphor-state-manager/greatlakes/sled-cycle.service
[Unit]
Description=Sled Cycle

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-state-manager/chassis-powercycle 0

[Install]
WantedBy=obmc-chassis-powercycle@0.target
