# meta-openembedded/meta-networking/recipes-support/tinyproxy/tinyproxy/tinyproxy.service
[Unit]
Description=Tinyproxy, Lightweight http(s) proxy daemon
RequiresMountsFor=/var

[Service]
Group=tinyproxy
ExecStartPre=/bin/mkdir -p /var/log/tinyproxy
ExecStartPre=/bin/chown -R nobody.tinyproxy /var/log/tinyproxy
ExecStartPre=/bin/mkdir -p /var/run/tinyproxy
ExecStartPre=/bin/chown -R nobody.tinyproxy /var/run/tinyproxy
ExecStart=/usr/bin/tinyproxy -d

[Install]
WantedBy=multi-user.target
