# meta-quanta/meta-olympus-nuvoton/recipes-x86/crashdump/files/com.intel.crashdump.service
[Unit]
Description=Intel BMC CPU Crashdump

[Service]
Restart=always
ExecStart={bindir}/crashdump
Type=simple

[Install]
WantedBy=multi-user.target
