# meta-facebook/meta-ventura/recipes-ventura/plat-svc/files/ventura-init-tray-sgpio-status.service
[Unit]
Description=Ventura Init Tray Sgpio Status
After=ventura-sys-init.service

[Service]
ExecStart=/usr/libexec/ventura/ventura-init-tray-sgpio-status
SyslogIdentifier=Ventura Init Tray Sgpio Status
Type=simple
Restart=on-failure
RestartSec=1

[Install]
WantedBy=multi-user.target
