# meta-facebook/meta-ventura/recipes-ventura/plat-svc/files/ventura-sys-init.service
[Unit]
Description=Ventura Early System Init
Before=phosphor-pid-control.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/ventura/ventura-early-sys-init
SyslogIdentifier=Ventura Early Init

[Install]
WantedBy=multi-user.target
