# meta-phosphor/common/recipes-phosphor/sensors/obmc-hwmon/obmc-hwmon.service
[Unit]
Description=Phosphor Hwmon Poller
Wants=mapper-wait@-org-openbmc-sensors.service
After=mapper-wait@-org-openbmc-sensors.service
Before=obmc-fans-ready.target

[Service]
Restart=always
ExecStart={sbindir}/hwmon.py
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
