# meta-phosphor/recipes-phosphor/interfaces/slpd-lite/slpd-lite.service
[Unit]
Description=Lightweight SLP Server

After=network.target

[Service]
ExecStart=/usr/bin/env slpd
SyslogIdentifier=slpd
Restart=always

[Install]
WantedBy=multi-user.target
