# meta-phosphor/recipes-phosphor/chassis/obmc-button-power/org.openbmc.buttons.Power@.service
[Unit]
Description=Phosphor Power Button%i

[Service]
Restart=always
ExecStart=/usr/bin/env button_power.exe
SyslogIdentifier=button_power.exe
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
