# meta-evb/meta-evb-nuvoton/meta-scm-npcm845/recipes-phosphor/fans/phosphor-pid-control/fan-boot-control.service
[Unit]
Description=Set Fan to Default Duty as Booting Up
DefaultDependencies=no

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/bin/fan-default-speed.sh

[Install]
WantedBy=sysinit.target
