# meta-evb/meta-evb-nuvoton/meta-scm-npcm845/recipes-phosphor/fans/phosphor-pid-control/phosphor-pid-control.service
[Unit]
Description=Fan Control Daemon
After=xyz.openbmc_project.fansensor.service

[Service]
Type=simple
ExecStart=/usr/bin/swampd
Restart=always
RestartSec=5
StartLimitInterval=0

[Install]
WantedBy=basic.target
