# meta-openbmc-machines/meta-openpower/meta-ibm/meta-romulus/recipes-phosphor/chassis/cpld-trigger/cpld_reset@.service
[Unit]
Description=Reset CPLD to wait after VDD is powered up
Wants=op-wait-power-off@%i.service
After=op-wait-power-off@%i.service
Conflicts=obmc-chassis-poweron@%i.target

[Service]
ExecStart=/usr/bin/env cpld_reset.sh
SyslogIdentifier=cpld_reset.sh
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=obmc-chassis-poweroff@%i.target
