# import-layers/meta-raspberrypi/recipes-connectivity/bluez5/bluez5/brcm43438.service
[Unit]
Description=Broadcom BCM43438 bluetooth HCI
ConditionPathIsDirectory=/proc/device-tree/soc/gpio@7e200000/bt_pins
Before=bluetooth.service
BindsTo=dev-serial1.device
After=dev-serial1.device

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/btuart

[Install]
WantedBy=multi-user.target
