# meta-phosphor/common/recipes-phosphor/obmc-phosphor-user/files/obmc-phosphor-user.service
[Unit]
Description=Phosphor OpenBMC user management daemon

[Service]
ExecStart=/usr/sbin/obmc-phosphor-userd
Environment="PYTHONUNBUFFERED=1"
Restart=always

[Install]
WantedBy=multi-user.target
