# import-layers/meta-openembedded/meta-networking/recipes-support/nis/ypbind-mt/ypbind.service
[Unit]
Description=YP Bind
Requires=domainname.service
After=domainname.service network.target
Before=systemd-user-sessions.service

[Service]
Type=forking
PIDFile=/var/run/ypbind.pid
ExecStart=/usr/sbin/ypbind

[Install]
WantedBy=multi-user.target

