# import-layers/meta-openembedded/meta-oe/recipes-navigation/gpsd/gpsd/gpsd.service
[Unit]
Description=GPS (Global Positioning System) Daemon
Requires=gpsd.socket

[Service]
EnvironmentFile=/etc/default/gpsd
ExecStart=/usr/sbin/gpsd -N $GPS_DEVICES

[Install]
Also=gpsd.socket
