# meta-facebook/meta-yosemite4/recipes-yosemite4/plat-tool/files/FanboardPowerOff@.service
[Unit]
Description=Power Off Fanboard %i
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/libexec/plat-tool/fanboardpower off %i
SyslogIdentifier=power-off-fanboard %i