# meta-ibm/recipes-phosphor/sensors/phosphor-hwmon/max31785-hwmon-helper@.service
[Unit]
Description=MAX31785 hwmon helper
Wants=phosphor-cooling-type@0.service
After=phosphor-cooling-type@0.service

[Service]
Type=oneshot
ExecStart=/usr/bin/env start_max31785_hwmon.sh %i
SyslogIdentifier=max31785-hwmon-helper
