# meta-openbmc-machines/meta-openpower/meta-ibm/recipes-phosphor/bmcweb/bmcweb/bmcweb.service
[Unit]
Description=BMC Redfish Server (bmcweb)
After=network.target
After=obmc-webserver-pre.target

[Service]
Restart=always
ExecStart=/usr/bin/env bmcweb
SyslogIdentifier=bmcweb
