# import-layers/meta-openembedded/meta-networking/recipes-protocols/quagga/files/isisd.service
[Unit]
Description=IS-IS routing daemon
BindTo=zebra.service
After=zebra.service
ConditionPathExists=@SYSCONFDIR@/quagga/isisd.conf

[Service]
Type=forking
EnvironmentFile=-@SYSCONFDIR@/default/quagga
ExecStart=@SBINDIR@/isisd -d $isisd_options -f /etc/quagga/isisd.conf
Restart=on-abort

[Install]
WantedBy=multi-user.target
