# meta-phosphor/recipes-support/rng-tools/rng-tools/rngd.service
[Unit]
Description=Hardware RNG Entropy Gatherer Daemon

[Service]
ExecStart=@SBINDIR@/rngd -f

[Install]
WantedBy=multi-user.target
