# meta-openembedded/meta-networking/recipes-protocols/net-snmp/net-snmp/snmpd.service
[Unit]
Description=Simple Network Management Protocol (SNMP) Daemon.
After=syslog.target network.target

[Service]
Type=notify
Environment=OPTIONS="-Ls0-6d"
EnvironmentFile=-/etc/default/snmpd
ExecStart=/usr/sbin/snmpd $OPTIONS -a -f
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
