# meta/recipes-connectivity/ppp/ppp/ppp@.service
[Unit]
Description=PPP link to %I
Before=network.target

[Service]
ExecStart=@SBINDIR@/pppd call %I nodetach nolog

[Install]
WantedBy=multi-user.target
