# meta-ibm/recipes-phosphor/fans/fan-watchdog/fan-watchdog.service
[Unit]
Description=Fan Watchdog Daemon
Wants=phosphor-fan-monitor-init@0.service
Before=phosphor-fan-monitor-init@0.service
Conflicts=obmc-fan-watchdog-takeover.target

[Service]
EnvironmentFile={envfiledir}/obmc/fan-watchdog/fan-watchdog.conf
ExecStart=/sbin/watchdog -T ${{TIMEOUT}} -t ${{INTERVAL}} -F ${{DEVICE}}
KillSignal=SIGKILL

[Install]
RequiredBy=obmc-chassis-poweron@0.target
