# meta-phosphor/common/recipes-phosphor/sensors/files/obmc-phosphor-sensord.service
[Unit]
Description=Phosphor OpenBMC sensor management daemon

[Service]
Restart=always
ExecStart=/usr/sbin/obmc-phosphor-sensord

[Install]
WantedBy=multi-user.target
