# meta-phosphor/common/recipes-phosphor/sensors/files/obmc-phosphor-sensord.service
[Unit]
Description=Phosphor Sensor Manager

[Service]
Restart=always
ExecStart={sbindir}/obmc-phosphor-sensord

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
