# meta-ampere/meta-common/recipes-ampere/host/ampere-hostctrl/ampere-host-on-host-check@.service
[Unit]
Description=Check Host%i status before obmc-host-startmin
After=obmc-host-started@%i.target
After=obmc-power-start-pre@%i.target
Conflicts=obmc-host-stop@%i.target
Conflicts=phosphor-reset-host-check@%i.service
OnFailure=obmc-host-quiesce@%i.target
OnSuccess=obmc-host-already-on@%i.target
OnFailureJobMode=flush

[Service]
Type=oneshot
ExecStart=/usr/sbin/ampere_host_check.sh
