# meta-ampere/meta-common/recipes-ampere/host/ampere-hostctrl/ampere-host-force-reset@.service
[Unit]
Description=Ampere Computing LLC force reset host%i
Wants=phosphor-reboot-host@%i.service
Conflicts=obmc-chassis-poweroff@%i.target
OnFailure=obmc-chassis-powerreset@%i.target

[Service]
RemainAfterExit=no
Type=oneshot
ExecStart=/bin/sh -c "sleep 3 && /usr/bin/env ampere_power_util.sh mb force_reset"
SyslogIdentifier=ampere_power_util.sh

[Install]
WantedBy=obmc-host-force-warm-reboot@%i.target