# meta-ampere/meta-common/recipes-ampere/host/ampere-hostctrl/obmc-power-already-on@.target
[Unit]
Description=Ampere power %i already on target
After=phosphor-wait-power-on@%i.service
Conflicts=obmc-chassis-poweroff@%i.target
Conflicts=ampere-host-shutdown-ack@%i.service
