# meta-phosphor/recipes-phosphor/inventory/phosphor-inventory-manager/xyz.openbmc_project.Inventory.Manager.service
[Unit]
Description=Phosphor Inventory Manager
Before=mapper-wait@-xyz-openbmc_project-inventory.service

[Service]
Restart=always
Type=dbus
ExecStart=/usr/bin/env phosphor-inventory
SyslogIdentifier=phosphor-inventory
BusName={BUSNAME}

[Install]
WantedBy=multi-user.target
