# meta-fii/meta-mori/recipes-mori/mori-fw-utility/mori-cpld-update/phosphor-ipmi-flash-mbcpld-update.service
[Unit]
Description=Phosphor-ipmi-flash update MB CPLD firmware service

[Service]
ExecStart=mori-fw.sh mbcpld /run/initramfs/mb-cpld.svf
Type=oneshot

[Install]
WantedBy=phosphor-ipmi-flash-mbcpld-update.target
