# meta-fii/meta-mori/recipes-mori/mori-fw-utility/mori-cpld-update/phosphor-ipmi-flash-bmccpld-update.service
[Unit]
Description=Phosphor-ipmi-flash update BMC CPLD firmware service

[Service]
ExecStart=mori-fw.sh bmccpld /run/initramfs/bmc-cpld.svf
Type=oneshot

[Install]
WantedBy=phosphor-ipmi-flash-bmccpld-update.target
