# meta-phosphor/common/recipes-phosphor/state/files/phosphor-reset-host-running@.service
[Unit]
Description=Start host%i running after BMC reset
After=obmc-host-reset-running@%i.target
Requires=obmc-host-reset-running@%i.target
ConditionPathExists=/run/openbmc/host@%i-on

[Service]
RemainAfterExit=no
ExecStart=/bin/systemctl start obmc-host-start@%i.target


[Install]
WantedBy=obmc-host-reset@%i.target
