# meta-phosphor/common/recipes-phosphor/state/files/xyz.openbmc_project.State.BMC.service
[Unit]
Description=Phosphor BMC State Manager
Wants=mapper-wait@-org-openbmc-managers-System.service
After=mapper-wait@-org-openbmc-managers-System.service
Before=mapper-wait@-xyz-openbmc_project-state-bmc.service

[Service]
ExecStart={sbindir}/phosphor-bmc-state-manager
Restart=always
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
