# meta-facebook/meta-harma/recipes-phosphor/state/phosphor-state-manager/chassis-poweron@.service
[Unit]
Description=power on chassis:%i

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-state-manager/chassis-poweron %i

[Install]
RequiredBy=obmc-chassis-poweron@%i.target
