# meta-openembedded/meta-networking/recipes-support/nis/yp-tools/domainname.service
[Unit]
Description=NIS Domainname

[Service]
Type=oneshot
EnvironmentFile=/etc/nisdomainname
ExecStart=/usr/bin/domainname $NISDOMAINNAME
RemainAfterExit=true

[Install]
WantedBy=multi-user.target

