# meta-facebook/meta-ventura2/recipes-ventura2/plat-svc/files/marvell-switch-init.service
[Unit]
Description=Marvell Switch Initialization
After=basic.target
Before=systemd-networkd.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/plat-svc/marvell-switch-init

[Install]
WantedBy=multi-user.target
