# meta-facebook/meta-ventura2/recipes-ventura2/plat-svc/files/sgpio-state-init.service
[Unit]
Description=SGPIO State Initialization

[Service]
Type=oneshot
ExecStart=/usr/libexec/plat-svc/sgpio-state-init
RemainAfterExit=yes
