# meta-facebook/meta-ventura2/recipes-ventura2/plat-svc/files/valve-dac-init@.service
[Unit]
Description=Valve DAC output voltage init (%i)

[Service]
Type=oneshot
ExecStart=/usr/libexec/plat-svc/valve-dac-init %i
RemainAfterExit=yes
