# meta-phosphor/recipes-phosphor/dbus/phosphor-objmgr/mapper-subtree-remove@.service
[Unit]
Description=mapper subtree-remove %I
Wants=obmc-mapper.target
After=obmc-mapper.target

[Service]
Type=oneshot
RemainAfterExit=no
Restart=on-failure
TimeoutStartSec=infinity
ExecStart=/usr/bin/mapper subtree-remove %I
