# import-layers/yocto-poky/meta/recipes-graphics/x11-common/xserver-nodm-init/xserver-nodm.service.in
[Unit]
Description=Xserver startup without a display manager

[Service]
EnvironmentFile=/etc/default/xserver-nodm
User=@USER@
ExecStart=/etc/xserver-nodm/Xserver

[Install]
Alias=display-manager.service
