# upstream-layers/openembedded-core/meta-selftest/recipes-test/cpp/files/cpp-example.service
[Unit]
Description=C++ Example Service
After=network.target

[Service]
Type=simple
User=@USER@
Group=@GROUP@
ExecStart=@BINDIR@/@BINARY_NAME@ --endless

[Install]
WantedBy=multi-user.target
