# meta-google/recipes-google/bare-metal-ipmi-kcs/bare-metal-ipmi-kcs/disable-ipmi-kcs.service.in
[Unit]
Description=Disable ipmi kcs
BindsTo=gbmc-bare-metal-active@0.target
Before=gbmc-bare-metal-active@0.target

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/bin/sh -c 'systemctl stop phosphor-ipmi-kcs@@KCS_DEV@'
ExecStop=/bin/sh -c 'systemctl start phosphor-ipmi-kcs@@KCS_DEV@'

[Install]
WantedBy=gbmc-bare-metal-active@0.target
