# meta-phosphor/common/recipes-phosphor/sensors/phosphor-hwmon/xyz.openbmc_project.Hwmon@.service
[Unit]
Description=Phosphor Hwmon Poller
ConditionFileNotEmpty=/etc/default/obmc/hwmon/%I.conf

[Service]
Restart=always
ExecStart={sbindir}/phosphor-hwmon-readd -o %f
EnvironmentFile=/etc/default/obmc/hwmon/%I.conf
