# meta-google/recipes-google/ncsi/files/gbmc-ncsi-sslh.socket.in
[Socket]
BindToDevice=@NCSI_IF@
ListenStream=3967

[Install]
WantedBy=sockets.target
