# meta-google/recipes-google/ncsi/files/gbmc-ncsi-old.service.in
[Unit]
Before=systemd-networkd.service
ConditionFileNotEmpty=/etc/systemd/network/00-bmc-@OLD_IF@.network

[Service]
Type=oneshot
ExecStart=mv -v /etc/systemd/network/00-bmc-@OLD_IF@.network /etc/systemd/network/00-bmc-@NCSI_IF@.network

[Install]
WantedBy=multi-user.target
