# meta-google/recipes-google/ncsi/files/gbmc-ncsi-ra.service.in
[Unit]
Description=RA Gateway Updater
Wants=gbmc-ncsi-networkd-wait.target
After=gbmc-ncsi-networkd-wait.target
After=ncsid@@NCSI_IF@.service
StartLimitIntervalSec=1min
StartLimitBurst=5

[Service]
Restart=always
RestartSec=5
ExecStart=/usr/libexec/gbmc-ncsi-ra.sh @NCSI_IF@

[Install]
WantedBy=multi-user.target
