# meta-google/recipes-google/ncsi/files/gbmc-ncsi-alias.service.in
[Unit]
BindsTo=ncsid@@NCSI_IF@.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/sbin/ip link set @NCSI_IF@ alias ncsi
ExecStop=/sbin/ip link set @NCSI_IF@ alias ''

[Install]
WantedBy=nic-hostful@@NCSI_IF@.target
WantedBy=nic-hostless@@NCSI_IF@.target
