# meta-google/recipes-google/ncsi/files/gbmc-ncsi-set-nicenabled.service.in
[Unit]
Description=Set NICEnabled property to true
Requires=gbmc-ncsi-networkd-wait.target
After=gbmc-ncsi-networkd-wait.target
StartLimitIntervalSec=10
StartLimitBurst=3

[Service]
Type=oneshot
ExecStart=busctl set-property xyz.openbmc_project.Network /xyz/openbmc_project/network/@NCSI_IF@ xyz.openbmc_project.Network.EthernetInterface NICEnabled b true
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
