# meta-ampere/meta-common/recipes-ac01/host/ampere-power-logger/service_files/xyz.openbmc_project.power_logger.service
[Unit]
Description=Host Monitor

[Service]
Restart=on-failure
ExecStart=/usr/bin/env ampere-power-logger
SyslogIdentifier=ampere-power-logger
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target