# meta-facebook/recipes-fb/fb-powerctrl/files/yosemitev2/host-gpio.service
[Unit]
Description=Configure GPIOs for Yosemitev2

[Service]
Restart=no
RemainAfterExit=true
Type=oneshot
ExecStart=/usr/sbin/setup_gpio.sh
SyslogIdentifier=setup_gpio.sh
