# meta-evb/meta-evb-nuvoton/meta-scm-npcm845/recipes-phosphor/gpio/cpu-err/cpu-err.service
[Unit]
Description=CPU ERR LED Service
Wants=mapper-wait@-xyz-openbmc_project-led-groups.service
After=mapper-wait@-xyz-openbmc_project-led-groups.service
StartLimitIntervalSec=0

[Service]
Type=oneshot
ExecStart=/usr/bin/toggle_caterr_led.sh
