# import-layers/meta-virtualization/recipes-extended/irqbalance/irqbalance/irqbalanced.service
[Unit]
Description=irqbalance daemon
After=syslog.target

[Service]
ExecStart=@SBINDIR@/irqbalance --foreground

[Install]
WantedBy=multi-user.target
