# meta-openbmc-machines/meta-x86/meta-intel/meta-common/recipes-intel/bmcweb/bmcweb/bmcweb.socket
[Unit]
Description=BMC Webserver socket

[Socket]
ListenStream=883
ReusePort=true

[Install]
WantedBy=sockets.target
