# meta-openembedded/meta-networking/recipes-support/open-vm-tools/open-vm-tools/vmtoolsd.service
[Unit]
Description=Open Virtual Machine Tools (VMware Tools)
ConditionVirtualization=vmware

[Service]
ExecStart=/usr/bin/vmtoolsd
Restart=on-failure
KillSignal=SIGKILL

[Install]
WantedBy=multi-user.target
