# meta-google/recipes-google/ssh/authorized-keys-comp/authorized-keys-comp.service
[Service]
Type=oneshot
ExecStart=/usr/libexec/authorized-keys-comp.sh

[Install]
WantedBy=multi-user.target
