# meta-quanta/meta-gbs/recipes-gbs/cpld-ver-check/files/cpld-version.service
[Unit]
Description=Report CPLD Version

[Service]
RemainAfterExit=yes
Type=oneshot
ExecStart=/usr/bin/cpld_version.sh

[Install]
WantedBy=multi-user.target
