# meta-phosphor/common/recipes-phosphor/examples/files/org.openbmc.examples.PythonService.service
[Unit]
Description=Phosphor Application Example

[Service]
Restart=always
ExecStart={sbindir}/obmc-phosphor-example-pydbus
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
