# upstream-layers/meta-openembedded/meta-networking/recipes-support/ntp/ntp/ntpd.service
[Unit]
Description=Network Time Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/sbin/ntpd -u ntp:ntp -n -g

[Install]
WantedBy=multi-user.target
