# upstream-layers/meta-openembedded/meta-oe/recipes-bsp/edac-utils/files/edac.service
[Unit]
Description=Initialize EDAC Drivers For Machine Hardware

[Service]
Type=oneshot
ExecStart=@SBINDIR@/edac-ctl --register-labels
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
