# meta-phosphor/common/recipes-phosphor/mount-overlay-help/mount-overlay-help/mount-machine-id.service
[Unit]
Description=Mount machine-id
Requires=etc.mount
After=etc.mount
DefaultDependencies=false
ConditionPathExists=!/var/persist/etc/machine-id

[Install]
RequiredBy=local-fs.target

[Service]
ExecStart=/bin/sh -c "mount --bind /run/machine-id /etc/machine-id"
Type=oneshot
RemainAfterExit=no
