# meta-phosphor/common/recipes-phosphor/mount-overlay-help/mount-overlay-help/prepare-overlay.service
[Unit]
Description=Prepare for overlay fs
RequiresMountsFor=/var
DefaultDependencies=false

[Service]
ExecStart=/bin/sh -c "if ! mountpoint /etc; then rm -rf /var/persist/etc-work && mkdir -p /var/persist/etc-work; fi; mkdir -p /var/persist/etc /var/persist/home/root"
Type=oneshot
RemainAfterExit=no
