# meta-ampere/meta-common/recipes-phosphor/flash/phosphor-software-manager/obmc-update-fru@.service
[Unit]
Description=Update image /tmp/images/%I/*.bin to FRU

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/sbin/ampere_firmware_upgrade.sh fru /tmp/images/%I/*.bin
