# meta-ampere/meta-common/recipes-phosphor/flash/phosphor-software-manager/obmc-flash-host-bios@.service
[Unit]
Description=Ampere service for flashing the Host firmware image

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/sbin/ampere_flash_bios.sh /tmp/images/%I/*.img
