# meta-ampere/meta-jade/recipes-phosphor/gpio/phosphor-gpio-monitor/ampere_scp_failover.service
[Unit]
Description=Ampere Computing LLC Ampere failover event service

[Service]
Restart=no
ExecStart=/usr/bin/env ampere_scp_failover.sh
