# meta-phosphor/recipes-phosphor/flash/obmc-mgr-download/org.openbmc.managers.Download.service
[Unit]
Description=Phosphor Download Manager

[Service]
Restart=always
ExecStart=/usr/bin/env download_manager.py
SyslogIdentifier=download_manager.py
Environment="PYTHONUNBUFFERED=1"
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy=multi-user.target
