# meta-axiado-sdk/recipes-axiado/tcu-reset/files/ax3000-reboot.service
[Unit]
Description=Use ax3000-fw-update for TCU reset
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/bin/tcu-reset.sh
RemainAfterExit=no

