# meta-openbmc-machines/meta-openpower/common/recipes-phosphor/host/op-pdbg-host-control/start_host@.service
[Unit]
Description=Start host%i SBE (pdbg)
After=obmc-host-start-pre@%i.target obmc-power-on@%i.target
Wants=obmc-host-start@%i.target obmc-power-on@%i.target
Before=obmc-host-start@%i.target obmc-host-started@%i.target
Conflicts=obmc-chassis-stop@%i.target

[Service]
ExecStart={bindir}/start_host.sh {MACHINE}
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=obmc-chassis-start@%i.target
