# meta-phosphor/common/recipes-phosphor/interfaces/phosphor-rocket/phosphor-rocket.service
[Unit]
Description=Phosphor Webserver
PartOf=systemd-networkd.service
After=obmc-webserver-pre.target

[Service]
Restart=always
ExecStart=/usr/bin/env phosphor-rocket $APPLICATION
SyslogIdentifier=phosphor-rocket
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile={envfiledir}/obmc/wsgi_app

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
