# yocto-poky/meta/recipes-kernel/oprofile/oprofileui-server/oprofileui-server.service
[Unit]
Description=OProfileUI Server
After=network.target

[Service]
ExecStart=/bin/sh -c ". @SYSCONFDIR@/profile; @BINDIR@/oprofile-server"
