# import-layers/meta-openembedded/meta-oe/recipes-extended/zram/zram/zram.service
[Unit]
Description=Enable zram compressed in-memory swap.
After=multi-user.target

[Service]
RemainAfterExit=yes
ExecStart=/usr/bin/zram-load.sh --load
ExecStop=/usr/bin/zram-load.sh --unload
Type=oneshot

[Install]
WantedBy=multi-user.target
