# meta-phosphor/common/recipes-phosphor/watchdog/phosphor-watchdog/phosphor-watchdog@.service
[Unit]
Description=Phosphor %I watchdog

[Service]
Restart=no
EnvironmentFile={envfiledir}/obmc/watchdog/%I
ExecStart={sbindir}/phosphor-watchdog --service=${{SERVICE}} --path=${{PATH}} --target=${{TARGET}}
