# meta-phosphor/common/recipes-phosphor/watchdog/phosphor-watchdog/phosphor-watchdog@.service
[Unit]
Description=Phosphor %I watchdog

[Service]
Restart=no
EnvironmentFile={envfiledir}/obmc/watchdog/%I
ExecStart=/usr/bin/env phosphor-watchdog --continue --service=${{SERVICE}} --path=${{DEVPATH}} --target=${{TARGET}}
SyslogIdentifier=phosphor-watchdog
