# upstream-layers/meta-openembedded/meta-networking/recipes-support/nuttcp/nuttcp/nuttcp.socket
[Unit]
Description=Network Performance Testing Server Activation Socket

[Socket]
ListenStream=5000
Accept=true

[Install]
WantedBy=sockets.target
