# upstream-layers/meta-openembedded/meta-networking/recipes-support/nuttcp/nuttcp/nuttcp@.service
[Unit]
Description=Network Performance Testing Server
After=local-fs.target

[Service]
User=nobody
ExecStart=/usr/bin/nuttcp -S
StandardInput=socket
