# meta-phosphor/common/recipes-phosphor/settings/settings/settings.service
[Unit]
Description=Phosphor Settings Manager

[Service]
ExecStart={sbindir}/settings_manager.py
Environment="PYTHONUNBUFFERED=1"
Restart=always

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
