# upstream-layers/meta-openembedded/meta-oe/recipes-extended/etcd/etcd/etcd-new.path
[Unit]
Description=Monitor the etcd config file changes

[Path]
PathChanged=/run/etcd-new.conf
Unit=etcd-new.service

[Install]
WantedBy=multi-user.target
