# meta-quanta/meta-olympus-nuvoton/recipes-support/nss-pam-ldapd/nss-pam-ldapd/nslcd.service
[Unit]
Description=LDAP daemon
After=syslog.target network.target
StartLimitIntervalSec=0

[Service]
Type=forking
ExecStart=/usr/sbin/nslcd

[Install]
WantedBy=multi-user.target
