# meta-phosphor/common/recipes-phosphor/chassis/obmc-watchdog/obmc-start-watchdog@.service
[Unit]
Description=Start Watchdog%i
Wants=obmc-host-started@%i.target
After=obmc-host-started@%i.target
Wants=mapper-wait@-org-openbmc-watchdog-host%i.service
After=mapper-wait@-org-openbmc-watchdog-host%i.service
Conflicts=obmc-chassis-stop@%i.target
ConditionPathExists=!/run/openbmc/host@%i-on

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart={sbindir}/mapper call /org/openbmc/watchdog/host%i org.openbmc.Watchdog set i 30000
ExecStart={sbindir}/mapper call /org/openbmc/watchdog/host%i org.openbmc.Watchdog start

[Install]
WantedBy=obmc-chassis-start@%i.target
