# meta-phosphor/recipes-phosphor/chassis/phosphor-skeleton-control-power/org.openbmc.control.Power@.service
[Unit]
Description=Phosphor Power%i Control

[Service]
EnvironmentFile={envfiledir}/obmc/power_control
Restart=always
ExecStart=/usr/bin/env power_control.exe $POLL_INTERVAL $PGOOD_TIMEOUT
SyslogIdentifier=power_control.exe
Type=dbus
BusName={BUSNAME}

[Install]
# Handled by bitbake recipe
