# meta-evb/meta-evb-nuvoton/meta-scm-npcm845/recipes-phosphor/console/obmc-console/obmc-console-ttyS1-ssh.socket
[Unit]
Description=Phosphor Host Console SSH Per-Connection socket
Wants=obmc-console@ttyS1.service

[Socket]
ListenStream=2220
Accept=yes

[Install]
WantedBy=sockets.target