# meta-fii/meta-kudo/recipes-phosphor/watchdog/phosphor-watchdog/phosphor-watchdog-host-reset.service
[Unit]
Description=Watchdog power state reset system

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=echo "Watchdog HardReset Action Happened (no ops)"
SyslogIdentifier=phosphor-watchdog
