# meta-fii/meta-kudo/recipes-phosphor/watchdog/phosphor-watchdog/phosphor-watchdog-host-cycle.service
[Unit]
Description=Watchdog power state power hotswap

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/sbin/kudo.sh rst hotswap
SyslogIdentifier=phosphor-watchdog
