# meta-fii/meta-kudo/recipes-phosphor/watchdog/phosphor-watchdog/phosphor-watchdog-host-reset.service
[Unit]
Description=Watchdog power state reset system

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/libexec/ampere-hostctrl/ampere_power_util.sh mb force_reset
SyslogIdentifier=phosphor-watchdog
