# meta-fii/meta-kudo/recipes-phosphor/watchdog/phosphor-watchdog/phosphor-watchdog-host-poweroff.service
[Unit]
Description=Watchdog power state power cycle shutdown

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/env kudo.sh rst shutdown
SyslogIdentifier=phosphor-watchdog
