# meta-facebook/recipes-phosphor/gpio/phosphor-gpio-monitor/obmc-ipmb-rescan-fru.service
[Unit]
Description=ipmb based rescan the slots fru service

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/ipmb-rescan-fru
