# meta-facebook/recipes-fb/reconfig-interface-duid-ll/reconfig-interface-duid-ll/reconfig-interface-duid-ll@.service
[Unit]
Description=Reconfig %i to get the IPv6 LL
Before=xyz.openbmc_project.Network.service

[Service]
ExecStart=/usr/libexec/reconfig-interface-duid-ll %i
SyslogIdentifier=reconfig-interface-duid-ll-%i
Type=oneshot

[Install]
WantedBy=multi-user.target
