# meta-openbmc-machines/meta-openpower/common/recipes-phosphor/vpd/openpower-fru-vpd/xyz.openbmc_project.OpenPOWER.Vpd.Parser@.service
[Unit]
Description=Read OpenPOWER-format VPD from EEPROM
Wants=mapper-wait@-xyz-openbmc_project-Inventory-Manager.service
After=mapper-wait@-xyz-openbmc_project-Inventory-Manager.service

[Service]
Restart=no
Type=oneshot
ExecStart={sbindir}/openpower-read-vpd --vpd=%f --fru=$FRUS --object=$PATHS
