# meta-fii/meta-kudo/recipes-phosphor/state/phosphor-state-manager/xyz.openbmc_project.State.Host.service
[Unit]
Description=Phosphor Host State Manager
Wants=mapper-wait@-xyz-openbmc_project-control-host0-auto_reboot.service
After=mapper-wait@-xyz-openbmc_project-control-host0-auto_reboot.service
Before=mapper-wait@-xyz-openbmc_project-state-host.service
Wants=obmc-mapper.target
After=obmc-mapper.target

[Service]
ExecStart=/usr/bin/phosphor-host-state-manager
Restart=always
Type=dbus
BusName=xyz.openbmc_project.State.Host

[Install]
WantedBy=multi-user.target
