# meta-facebook/meta-ventura/recipes-phosphor/sensors/dbus-sensors/xyz.openbmc_project.leakdetector.critical@.service
[Unit]
Description=%i Critical Leak Handler

[Service]
Type=oneshot
ExecStart=/usr/libexec/dbus-sensors/critical-leak-handler %i
SyslogIdentifier=%i-critical-leak-handler
