# meta-facebook/meta-ventura/recipes-phosphor/sensors/dbus-sensors/xyz.openbmc_project.leakdetector.critical.deassert@.service
[Unit]
Description=%i Critical deassert Handler

[Service]
Type=oneshot
ExecStart=/usr/libexec/dbus-sensors/deassert-leak-handler %i
SyslogIdentifier=%i-deaasert-leak-handler
