# meta-quanta/meta-olympus-nuvoton/recipes-phosphor/sensors/phosphor-hwmon/olympus-reload-sensor-on.service
[Unit]
Description=Reload Sensor On Configuration

[Service]
ExecStart=/bin/sh -c "olympus-reload-sensor.sh 1"
SyslogIdentifier=obmc-power-start
