# meta-quanta/meta-olympus-nuvoton/recipes-phosphor/sensors/phosphor-hwmon/olympus-reload-sensor-off.service
[Unit]
Description=Reload Sensor off Configuration

[Service]
ExecStart=/bin/sh -c "olympus-reload-sensor.sh 0"
SyslogIdentifier=obmc-power-stop
