# meta-facebook/meta-tiogapass/recipes-tiogapass/fb-powerctrl/files/host-poweron.service
[Unit]
Description=Poweron Host Server
Requires=host-gpio.service
After=host-gpio.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/power-util mb on
SyslogIdentifier=power-util

[Install]
WantedBy=basic.target
