# meta-openbmc-machines/meta-x86/meta-quanta/meta-q71l/recipes-phosphor/quanta-powerctrl/files/host-gpio.service
[Unit]
Description=Configure GPIOs for Host Power Control

[Service]
Restart=no
RemainsAfterExit=true
Type=oneshot
ExecStart={sbindir}/init_once.sh
