# meta-openbmc-machines/meta-x86/meta-quanta/meta-q71l/recipes-phosphor/quanta-powerctrl/files/host-poweroff.service
[Unit]
Description=Stop Host
Requires=host-gpio.service
After=host-gpio.service
Conflicts=obmc-host-start@0.target

[Service]
RemainAfterExit=yes
Type=oneshot
ExecStart={sbindir}/poweroff.sh

[Install]
WantedBy=obmc-host-stop@0.target
