# meta-openbmc-machines/meta-x86/meta-quanta/meta-q71l/recipes-phosphor/quanta-powerctrl/files/host-poweroff.service
[Unit]
Description=Stop Host
Requires=host-gpio.service
After=host-gpio.service

[Service]
RemainAfterExit=yes
Type=oneshot
ExecStart=/usr/sbin/poweroff.sh

[Install]
WantedBy=obmc-chassis-stop@0.target
