# meta-facebook/meta-catalina/recipes-phosphor/gpio/phosphor-gpio-monitor/assert-gpio-log@.service
[Unit]
Description=%i GPIO Assert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/catalina-gpio-monitor %i assert

SyslogIdentifier=%i-gpio-assert
