# meta-facebook/meta-catalina/recipes-phosphor/gpio/phosphor-gpio-monitor/assert-uart-select-led.service
[Unit]
Description=Assert Uart Select

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/set-uart-select-led 1

SyslogIdentifier=assert-uart-select
