# meta-facebook/meta-catalina/recipes-phosphor/gpio/phosphor-gpio-monitor/assert-reset-button.service
[Unit]
Description=Assert ID_RST_BTN_BMC_N

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/catalina-gpio-monitor ID_RST_BTN_BMC_N assert

SyslogIdentifier=assert-run-power-pg
