# meta-facebook/meta-catalina/recipes-phosphor/gpio/phosphor-gpio-monitor/assert-run-power-pg.service
[Unit]
Description=Assert RUN_POWER_PG

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/catalina-gpio-monitor RUN_POWER_PG assert

SyslogIdentifier=assert-run-power-pg
