# meta-facebook/meta-catalina/recipes-phosphor/gpio/phosphor-gpio-monitor/assert-cpu-boot-done.service
[Unit]
Description=Assert CPU_BOOT_DONE

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/catalina-gpio-monitor CPU_BOOT_DONE assert

SyslogIdentifier=assert-cpu-boot-done