# meta-facebook/meta-catalina/recipes-phosphor/gpio/phosphor-gpio-monitor/catalina-host-ready.target
[Unit]
Description=Catalina host ready target
