# meta-facebook/meta-catalina/recipes-phosphor/gpio/phosphor-gpio-monitor/deassert-reset-button.service
[Unit]
Description=Deassert ID_RST_BTN_BMC_N

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/catalina-gpio-monitor ID_RST_BTN_BMC_N deassert

SyslogIdentifier=deassert-run-power-pg
