# meta-facebook/meta-catalina/recipes-phosphor/gpio/phosphor-gpio-monitor/deassert-uart-select-led.service
[Unit]
Description=Deassert Uart Select

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/set-uart-select-led 0

SyslogIdentifier=deassert-uart-select
