# meta-facebook/meta-catalina/recipes-phosphor/gpio/phosphor-gpio-monitor/platform-host-ready.target
[Unit]
Description=Catalina/Clemente host ready target
