# meta-facebook/meta-catalina/recipes-phosphor/gpio/phosphor-gpio-monitor/vr-fault-assert-log@.service
[Unit]
Description=%i VR Fault Assert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/vr-fault-event-logger "assert" "%i"
SyslogIdentifier=%i-vr-fault-assert
