# meta-facebook/meta-santabarbara/recipes-santabarbara/plat-svc/files/santabarbara-eid-init.service
[Unit]
Description=Set I2C EID at boot
DefaultDependencies=no
Before=xyz.openbmc_project.EntityManager.service
Before=mctp_init.service
Before=pldm.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/santabarbara-eid-init
SyslogIdentifier=Santabarbara Eid Init

[Install]
WantedBy=multi-user.target
