# import-layers/meta-openembedded/meta-oe/recipes-devtools/cloud9/cloud9/cloud9.service
[Unit]
Description=Cloud9 IDE
ConditionPathExists=|/var/lib/cloud9

[Service]
Restart=always
EnvironmentFile=-/etc/default/node
ExecStart=/usr/bin/node4 /usr/share/cloud9/bin/cloud9.js -l 0.0.0.0 -w /var/lib/cloud9 -p 3000

[Install]
WantedBy=multi-user.target
