# meta-phosphor/common/recipes-phosphor/host-ipmid/host-ipmi-hw-example/host-ipmi-hw-example.service
[Unit]
Description=Phosphor OpenBMC host IPMI to DBUS example implementation

[Service]
Restart=always
ExecStart=/usr/sbin/host-ipmi-hw

[Install]
WantedBy=multi-user.target
