# meta-quanta/meta-olympus-nuvoton/recipes-olympus-nuvoton/ipmi-bmc-update/ipmi-bmc-update/phosphor-ipmi-flash-bmc-verify.service
[Unit]
Description=Phosphor-ipmi-flash bmc verify service

[Service]
ExecStart=/usr/bin/bmc-verify.sh
Type=oneshot

[Install]
WantedBy=phosphor-ipmi-flash-bmc-verify.target
