# meta-evb/meta-evb-nuvoton/meta-buv-runbmc/recipes-phosphor/gpio/phosphor-gpio-monitor/phosphor-multi-gpio-monitor.service
[Unit]
Description=Phosphor GPIO monitor
After=xyz.openbmc_project.EntityManager.service

[Service]
Type=simple
Restart=always
RestartSec=5
StartLimitInterval=0
ExecStart=/usr/bin/phosphor-multi-gpio-monitor --config /usr/share/phosphor-gpio-monitor/BUV-GpioMonitorConfig.json

[Install]
WantedBy=multi-user.target
