# meta-evb/meta-evb-nuvoton/meta-buv-runbmc/recipes-phosphor/gpio/phosphor-gpio-monitor/psusensor_reload.service
[Unit]
Description=psusensor reload

[Service]
Type=oneshot
ExecStart=/bin/sh -c "psu_hotplug_action.sh 1"
SyslogIdentifier=psusensor_reload
