# meta-ampere/meta-common/recipes-ac01/host/ampere-host-state/xyz.openbmc_project.AmpHostState.service
[Unit]
Description=Ampere Host State verifying and updating
Wants=xyz.openbmc_project.State.Host.service
After=xyz.openbmc_project.State.Host.service

[Service]
Restart=always
ExecStart=/usr/bin/ampere-host-state
Type=simple

[Install]
WantedBy=multi-user.target
