# meta-openpower/recipes-phosphor/flash/openpower-software-manager/obmc-flash-bios-ubiattach.service
[Unit]
Description=Attach the PNOR mtd device to UBI

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/env obmc-flash-bios ubiattach
SyslogIdentifier=obmc-flash-bios
