# meta-phosphor/common/recipes-phosphor/obmc-phosphor-fan/files/obmc-phosphor-fand.service
[Unit]
Description=Phosphor OpenBMC fan management daemon

[Service]
Restart=always
ExecStart=/usr/sbin/obmc-phosphor-fand

[Install]
WantedBy=multi-user.target
