# meta-facebook/meta-yosemite4/recipes-phosphor/fans/phosphor-pid-control/monitor-pldm-sensor.service
[Unit]
Description=Monitor PLDM Sensors
After=xyz.openbmc_project.FruDevice.service xyz.openbmc_project.PLDM.service
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-pid-control/monitor-pldm-sensor

[Install]
WantedBy=multi-user.target
