# meta-phosphor/common/recipes-phosphor/interfaces/rest-dbus/rest-dbus.service
[Unit]
Description=Phosphor DBUS Browser

[Service]
Restart=always
ExecStart=/usr/bin/env rest-dbus
SyslogIdentifier=rest-dbus
Environment="PYTHONUNBUFFERED=1"
