# meta-phosphor/common/recipes-phosphor/interfaces/rest-dbus/rest-dbus.socket
[Unit]
Description=Phosphor DBUS Browser socket

[Socket]
ListenStream=3000

[Install]
WantedBy=sockets.target
