# meta-quanta/meta-olympus-nuvoton/recipes-phosphor/fans/phosphor-pid-control/phosphor-pid-control-stop.service
[Unit]
Description=OpenBMC Fan Control Daemon (Stop)

[Service]
Type=oneshot
ExecStart=-/bin/systemctl stop phosphor-pid-control.service
