# meta-quanta/meta-olympus-nuvoton/recipes-olympus-nuvoton/network/olympus-nuvoton-iptable-restore/phosphor-restore-iptable-rules.service
[Unit]
Description=Restore iptable rules
DefaultDependencies=no
ConditionPathExists=/var/lib/iptables_rules

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/bin/sh -c "iptables-restore < /var/lib/iptables_rules"

[Install]
WantedBy=multi-user.target
