# meta-openbmc-machines/meta-openpower/meta-ingrasys/meta-zaius/recipes-phosphor/chassis/avsbus-control/avsbus-disable@.service
[Unit]
Description=Disable AVSBus on VRMs
Wants=obmc-power-stop-pre@%i.target
Before=obmc-power-stop-pre@%i.target
Conflicts=obmc-chassis-start@%i.target
ConditionPathExists=!/run/openbmc/chassis@%i-on

[Service]
ExecStart={bindir}/zaius_avsbus.sh disable
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=obmc-chassis-stop@%i.target
