# meta-axiado-sdk/recipes-axiado/ax-usb/ax-usb/create-usb-rndis.service
[Unit]
Description=Setup USB RNDIS gadget function
Requires=bmc-boot-complete.service
After=bmc-boot-complete.service

[Service]
Type=oneshot
ExecStart=/usr/bin/create-usb-rndis.sh

[Install]
WantedBy=multi-user.target
