# meta-hpe/meta-rl300-g11/recipes-hpe/power-supplies-monitoring/files/psu1down.service
[Unit]
Description=HPE RL300 PSU1 Manager

[Service]
ExecStart=/usr/bin/psus-manager.sh 1 down

[Install]
# WantedBy=multi-user.target
