# meta-hpe/meta-rl300-g11/recipes-hpe/power-supplies-monitoring/files/psu2down.service
[Unit]
Description=HPE RL300 PSU2 Manager

[Service]
ExecStart=/usr/bin/psus-manager.sh 2 down

[Install]
# WantedBy=multi-user.target
