# meta-phosphor/common/recipes-phosphor/chassis/obmc-op-control-power/op-reset-chassis-on@.service
[Unit]
Description=Start chassis%i on after BMC reset
Requires=op-reset-set-power-on@%i.service
After=op-reset-set-power-on@%i.service
After=obmc-power-reset-on@%i.target
Requires=obmc-power-reset-on@%i.target

[Service]
RemainAfterExit=no
ExecStart=/bin/systemctl start obmc-power-chassis-on@%i.target


[Install]
WantedBy=obmc-chassis-reset@%i.target
