# meta-phosphor/common/recipes-phosphor/chassis/obmc-op-control-power/op-reset-chassis-on@.service
[Unit]
Description=Start chassis%i on after BMC reset
Requires=op-reset-chassis-running@%i.service
After=op-reset-chassis-running@%i.service
After=obmc-power-reset-on@%i.target
Requires=obmc-power-reset-on@%i.target
ConditionPathExists=/run/openbmc/chassis@%i-on

[Service]
RemainAfterExit=no
ExecStart=/bin/systemctl start obmc-chassis-poweron@%i.target


[Install]
WantedBy=obmc-chassis-powerreset@%i.target
