# meta-phosphor/common/recipes-phosphor/chassis/obmc-op-control-power/op-reset-set-power-on@.service
[Unit]
Description=Create file to indicate Pgood%i is on
After=obmc-chassis-reset-pgood-on@%i.target

[Service]
RemainAfterExit=no
Type=oneshot
ExecStart=/bin/sh -c "mkdir -p /run/openbmc/ && touch /run/openbmc/chassis@%i-on"


[Install]
WantedBy=obmc-chassis-reset@%i.target
