# meta-phosphor/common/recipes-phosphor/chassis/obmc-op-control-power/op-reset-set-power-on@.service
[Unit]
Description=Create file to indicate chassis%i pgood is on
After=obmc-power-reset-on@%i.target
Requires=obmc-power-reset-on@%i.target
After=op-reset-pgood-check@%i.service
Requires=op-reset-pgood-check@%i.service
Conflicts=obmc-power-chassis-off@%i.target

[Service]
RemainAfterExit=no
Type=oneshot
ExecStart=/bin/sh -c "mkdir -p /run/openbmc/ && touch /run/openbmc/chassis@%i-on"


[Install]
WantedBy=obmc-chassis-reset@%i.target
