# meta-facebook/meta-yosemite5/recipes-phosphor/dbus/fan-policy/fan-fault-policy.service
[Unit]
Description=Fan fault handler service for host protection

[Service]
Type=oneshot
ExecStart=/usr/libexec/fan-policy/fan-fault-policy