# meta-evb/meta-evb-nuvoton/meta-scm-npcm845/recipes-scm-npcm845/reload-sensors/files/reload-sensors.service
[Unit]
Description=Reload sensors

[Service]
ExecStart=/usr/bin/reload-sensors.sh
StandardOutput=syslog
Type=simple
Restart=on-failure

[Install]
WantedBy=obmc-post-complete.target
