# meta-openpower/recipes-phosphor/host/aspeed-lpc-ctrl/pnorboot.service
[Unit]
Description=ASPEED LPC boot from PNOR

[Service]
Type=oneshot
ExecStart=/usr/sbin/pnorboot --size 32

[Install]
WantedBy=multi-user.target
