# meta-phosphor/recipes-phosphor/interfaces/bmcweb/bmcweb.socket
[Unit]
Description=BMC Webserver socket

[Socket]
ListenStream=443
ReusePort=true

[Install]
WantedBy=sockets.target
