# meta-ampere/meta-mitchell/recipes-phosphor/fans/phosphor-fan/phosphor-fan-monitor@.service
[Unit]
Description=Phosphor Fan Monitor Daemon
Wants=xyz.openbmc_project.fansensor.service
After=xyz.openbmc_project.fansensor.service
StartLimitIntervalSec=30
StartLimitBurst=3

[Service]
Restart=always
ExecStart=-phosphor-fan-monitor
RestartSec=20
