# meta-ampere/meta-jade/recipes-ampere/platform/ampere-mac-update/ampere_update_mac.service
[Unit]
Description=Updating MAC Address Service From FRU inventory
Requires=xyz.openbmc_project.EntityManager.service
After=xyz.openbmc_project.EntityManager.service

[Service]
ExecStart=/usr/sbin/ampere_update_mac.sh eth1 3 80
SyslogIdentifier=ampere_update_mac.sh
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
