# meta-facebook/meta-yosemite5/recipes-phosphor/state/phosphor-state-manager/chassis-poweroff@.service
[Unit]
Description=Chassis Power Off: %i

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-state-manager/chassis-poweroff %i

[Install]
WantedBy=obmc-chassis-poweroff@%i.target
