# meta-facebook/meta-yosemite5/recipes-phosphor/state/phosphor-state-manager/chassis-poweron@.service
[Unit]
Description=Chassis Power On: %i

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-state-manager/chassis-poweron %i

[Install]
WantedBy=obmc-chassis-poweron@%i.target
