# meta-facebook/meta-yosemite5/recipes-phosphor/state/phosphor-state-manager/host-graceful-poweroff-failure@.service
[Unit]
Description=Handle host graceful poweroff failure for instance %i

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-state-manager/host-graceful-poweroff-failure %i

[Install]
WantedBy=multi-user.target