# meta-facebook/meta-yosemite5/recipes-phosphor/state/phosphor-state-manager/host-graceful-poweroff@.service
[Unit]
Description=power off host:%i
Wants=obmc-host-stop-pre@%i.target
Before=obmc-host-stop-pre@%i.target
Conflicts=obmc-host-start@%i.target
ConditionPathExists=!/run/openbmc/host@%i-request
OnFailure=host-graceful-poweroff-failure@%i.service

[Service]
Restart=no
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/phosphor-state-manager/host-graceful-poweroff %i

[Install]
WantedBy=obmc-host-warm-reboot@%i.target
