# meta-phosphor/recipes-phosphor/settings/phosphor-settings-manager/xyz.openbmc_project.Settings.service
[Unit]
Description=Phosphor Settings Daemon

[Service]
ExecStart=/usr/bin/env phosphor-settings-manager
SyslogIdentifier=phosphor-settings-manager
Restart=always
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy=multi-user.target
