# meta-facebook/meta-yosemite5/recipes-phosphor/dbus/thermal-policy/thermal-unr-policy.service
[Unit]
Description=Thermal UNR handler service for host protection

[Service]
Type=oneshot
ExecStart=/usr/libexec/thermal-policy/thermal-unr-policy