# meta-phosphor/common/recipes-phosphor/phosphor-led-manager/phosphor-led-manager/xyz.openbmc_project.ledmanager.service
[Unit]
Description=Phosphor LED Manager daemon

[Service]
Restart=always
ExecStart={sbindir}/phosphor-ledmanager
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
