# meta-facebook/recipes-phosphor/fans/phosphor-fan/obmc-chassis-poweroff@0.target
[Unit]
Description=Chassis%i (Power Off)
After=multi-user.target
