# meta-facebook/recipes-phosphor/fans/phosphor-fan/obmc-chassis-hard-poweroff@.target
[Unit]
Description=Chassis%i (Hard Power Off)
After=multi-user.target
Requires=obmc-poweroff.service
