# meta-facebook/recipes-phosphor/fans/phosphor-fan/obmc-poweroff.service
[Unit]
Description=power off host

[Service]
RemainAfterExit=no
Type=oneshot
ExecStart=/usr/libexec/phosphor-fan-sensor-monitor/host-poweroff
