# meta-fii/meta-mori/recipes-mori/mori-fw-utility/mori-bios-update/phosphor-ipmi-flash-bios-update.service
[Unit]
Description=Phosphor-ipmi-flash update BIOS firmware service

[Service]
ExecStart=mori-fw.sh bios /run/initramfs/bios-image
Type=oneshot
