# meta-openembedded/meta-networking/recipes-daemons/dnrd/dnrd/dnrd.service
[Unit]
Description=Domain Name Relay Daemon
After=network.target

[Service]
Type=forking
PIDFile=/var/run/dnrd.pid
EnvironmentFile=/etc/dnrd/dnrd.conf
ExecStart=/usr/sbin/dnrd $DNRD_OPTS -u $DNRD_USER

[Install]
WantedBy=multi-user.target
