# meta-phosphor/common/recipes-phosphor/chassis/obmc-mgr-inventory/obmc-sync-inventory-item@.service
[Unit]
Description=Sync %I
Wants=obmc-read-eeprom@%i.service
After=obmc-read-eeprom@%i.service

[Service]
Restart=no
Type=oneshot

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
