# meta-phosphor/recipes-phosphor/trace-enable/trace-enable/trace-enable.service
[Unit]
Description=Enable Linux trace events in the boot loader
RequiresMountsFor=/etc

[Service]
Type=oneshot
ExecStart=/usr/libexec/trace-enable /etc/trace-events.conf

[Install]
WantedBy=multi-user.target
