# meta-inventec/meta-transformers/recipes-phosphor/init/transformers-init/transformers-init.service
[Unit]
Description=Essential init commands for transformers

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/sbin/transformers-init.sh

[Install]
WantedBy=multi-user.target
