# meta-phosphor/recipes-devtools/iotools/iotools/iotools-setup.service
[Unit]
Description=Install iotool links

[Service]
RemainAfterExit=no
Type=oneshot
ExecStart=/usr/bin/env iotools --make-links
SyslogIdentifier=iotools

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
