# meta-facebook/meta-minerva/recipes-minerva/plat-svc/files/minerva-sys-init.service
[Unit]
Description=Minerva Early System Init
Before=phosphor-multi-gpio-monitor.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/minerva-early-sys-init
SyslogIdentifier=Minerva Early Init

[Install]
WantedBy=multi-user.target
