# import-layers/yocto-poky/meta/recipes-connectivity/portmap/portmap/portmap.service
[Unit]
Description=The RPC portmapper
After=network.target

[Service]
Type=forking
ExecStart=@BASE_SBINDIR@/portmap

[Install]
WantedBy=multi-user.target
