# upstream-layers/meta-openembedded/meta-filesystems/dynamic-layers/meta-python/recipes-support/gpiod-sysfs-proxy/gpiod-sysfs-proxy/run-gpio-sys.mount
# SPDX-License-Identifier: CC0-1.0
# SPDX-FileCopyrightText: 2024 Bartosz Golaszewski <bartosz.golaszewski@linaro.org>

[Unit]
Description=Remount of sysfs for gpiod-sysfs-proxy
ConditionPathExists=!/sys/class/gpio

[Mount]
DirectoryMode=0700
What=sysfs
Where=/run/gpio/sys
Type=sysfs
Options=nosuid,nodev,noexec
