# import-layers/meta-openembedded/meta-python/recipes-connectivity/gateone/gateone/gateone.service.in
[Unit]
Description=GateOne daemon
ConditionPathExists=|@localstate@/lib/gateone

[Service]
WorkingDirectory=@localstate@/lib/gateone
ExecStart=@bindir@/python @bindir@/gateone

[Install]
WantedBy=multi-user.target
