# meta-evb/meta-evb-nuvoton/meta-scm-npcm845/recipes-phosphor/ipmi/phosphor-ipmi-host/phosphor-ipmi-warm-reset.target
[Unit]
Description=IPMI service collection
AllowIsolate=no
