# meta-phosphor/common/recipes-phosphor/dbus/obmc-rest/obmc-rest.service
[Unit]
Description=Phosphor OpenBMC DBus REST daemon
Requires=obmc-mapper.service
After=obmc-mapper.service
PartOf=systemd-networkd.service

[Service]
Restart=always
ExecStart=/usr/sbin/obmc-rest

[Install]
WantedBy=multi-user.target
