# meta-phosphor/common/recipes-phosphor/dbus/obmc-rest/obmc-rest.socket
[Unit]
Description=Phosphor OpenBMC DBus REST daemon socket

[Socket]
ListenStream=443

[Install]
WantedBy=sockets.target
