# meta-phosphor/common/recipes-core/systemd/systemd/obmc-standby-reached.target
[Unit]
Description=OpenBMC has reached standby
After=obmc-standby.target
