# meta-phosphor/common/recipes-core/systemd/systemd/obmc-mapper.target
[Unit]
Description=OpenBMC Object Mapper
RefuseManualStart=yes
