# meta-phosphor/common/recipes-core/systemd/systemd/obmc-standby.target
[Unit]
Description=OpenBMC System
RefuseManualStart=yes
