# meta-fii/meta-kudo/recipes-kudo/kudo-fw-utility/kudo-fw/kudo-fw-ver.service
[Unit]
Description=Firmware Version Daemon

[Service]
RemainAfterExit=yes
Type=oneshot
ExecStart=/usr/libexec/kudo-fw/kudo-fw-ver.sh

[Install]
WantedBy=multi-user.target
