# meta-phosphor/recipes-phosphor/users/phosphor-user-manager/xyz.openbmc_project.LDAP.PrivilegeMapper.service
[Unit]
Description=Phosphor LDAP privilege mapper

[Service]
ExecStart=/usr/bin/env phosphor-ldap-mapper
SyslogIdentifier=phosphor-ldap-mapper
Restart=always
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy=multi-user.target
