# meta-phosphor/common/recipes-phosphor/net-ipmid/phosphor-net-ipmid/phosphor-net-ipmid.socket
[Socket]
ListenDatagram=623

[Install]
WantedBy=sockets.target

