# meta-arm/meta-arm/recipes-security/optee/optee-client/tee-supplicant.service
[Unit]
Description=TEE Supplicant
ConditionPathExistsGlob=/dev/teepriv[0-9]*

[Service]
User=root
EnvironmentFile=-@sysconfdir@/default/tee-supplicant
ExecStart=@sbindir@/tee-supplicant $OPTARGS

[Install]
WantedBy=basic.target
