# poky/meta/recipes-connectivity/openssh/openssh/sshdgenkeys.service
[Unit]
Description=OpenSSH Key Generation
RequiresMountsFor=/var /run

[Service]
ExecStart=@LIBEXECDIR@/sshd_check_keys
Type=oneshot
RemainAfterExit=yes
Nice=10
