# meta-phosphor/common/recipes-connectivity/avahi/avahi/check-avahi-pre-cond.service
[Unit]
Description=Check the precondition before starting avahi
Before=avahi-daemon.service

[Service]
RemainAfterExit=no
Type=oneshot
ExecStart=/usr/sbin/check-local-domain
ExecStart=/usr/sbin/create-user
ExecStart=/bin/touch /tmp/avahi

[Install]
WantedBy=obmc-standby.target
