# meta-phosphor/common/recipes-phosphor/obmc-phosphor-time-manager/files/obmc-phosphor-time-manager.service
[Unit]
Description=Phosphor Time Manager daemon
Requires=network.service
After=network.service
Requires=settings.service
After=settings.service

[Service]
Restart=always
ExecStart=/usr/sbin/timemanager

[Install]
WantedBy=multi-user.target
