# meta-facebook/meta-sanmiguel/recipes-phosphor/state/phosphor-state-manager/host-powerreset@.service
[Unit]
Description=Host Power Reset: %i

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-state-manager/host-powerreset %i

[Install]
WantedBy=obmc-host-force-warm-reboot@%i.target
