# meta-evb/meta-evb-nuvoton/meta-scm-npcm845/recipes-scm-npcm845/wdat/files/wdat@.service
[Unit]
Description=WDAT service

[Service]
Type=simple
StandardOutput=journal
SyslogIdentifier=wdat
ExecStart=/usr/sbin/wdat /dev/raw-$i
