# meta-openbmc-machines/meta-openpower/common/recipes-phosphor/occ/openpower-occ-control/op-occ-disable@.service
[Unit]
Description=OpenPOWER OCC Active Disable
Before=obmc-host-stop-pre@%i.target
Conflicts=op-occ-enable@%i.service

[Service]
ExecStart=/usr/bin/env occ-active.sh disable
SyslogIdentifier=occ-active.sh
Type=oneshot

[Install]
WantedBy=obmc-host-stop@%i.target
