# meta-phosphor/common/recipes-phosphor/flash/phosphor-software-manager/obmc-flash-bmc-ubiro@.service
[Unit]
Description=Create read-only ubi volumes %I

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart={sbindir}/obmc-flash-bmc ubiro {RO_MTD} ro-%i %i
ExecStart={sbindir}/obmc-flash-bmc ubikernel {KERNEL_MTD} kernel-%i %i
