# meta-phosphor/common/recipes-phosphor/flash/phosphor-software-manager/obmc-flash-bmc-ubiro-remove@.service
[Unit]
Description=Deletes read-only and kernel ubi volume %I

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart={sbindir}/obmc-flash-bmc ubiremove ro-%i
ExecStart={sbindir}/obmc-flash-bmc ubiremove kernel-%i
