# meta-phosphor/common/recipes-phosphor/flash/phosphor-software-manager/obmc-flash-bmc-ubirw-reset.service
[Unit]
Description=Reset BMC read-write ubi volume

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/bin/sh -c 'if [ $(fw_printenv rwreset) == "rwreset=true" ]; then \
                          obmc-flash-bmc ubirw {RW_MTD} rw; fi'
