# meta-phosphor/common/recipes-phosphor/flash/phosphor-software-manager/obmc-flash-bmc-ubirw.service
[Unit]
Description=Create BMC read-write ubi volume

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart={sbindir}/obmc-flash-bmc ubirw {RW_MTD} rw
