# meta-phosphor/common/recipes-phosphor/flash/phosphor-software-manager/usr-local.automount
[Unit]
Description=Automount /usr/local
After=obmc-standby.target

[Automount]
Where=/usr/local

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
