# meta-phosphor/common/recipes-phosphor/flash/phosphor-software-manager/usr-local.mount
[Unit]
Description=/usr/local
After=obmc-standby.target

[Mount]
What=tmpfs
Where=/usr/local
Type=tmpfs
