# meta-phosphor/common/recipes-phosphor/flash/phosphor-software-manager/xyz.openbmc_project.Software.Download.service
[Unit]
Description=Phosphor Download Manager

[Service]
ExecStart={sbindir}/phosphor-download-manager
Restart=always
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
