# meta-phosphor/common/recipes-phosphor/flash/phosphor-software-manager/xyz.openbmc_project.Software.Version.service
[Unit]
Description=Phosphor Version Software Manager
Wants=xyz.openbmc_project.Software.BMC.Updater.service
Before=xyz.openbmc_project.Software.BMC.Updater.service

[Service]
ExecStartPre=/bin/sh -c 'mkdir -p /tmp/images'
ExecStart={sbindir}/phosphor-version-software-manager
Restart=always
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
