# meta-phosphor/common/recipes-phosphor/flash/phosphor-software-manager/xyz.openbmc_project.Software.Version.service
[Unit]
Description=Phosphor Version Software Manager
Before=mapper-wait@-xyz-openbmc_project-software-version.service

[Service]
ExecStart={sbindir}/phosphor-version-software-manager
Restart=always
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
