# meta-phosphor/common/recipes-phosphor/flash/phosphor-software-manager/xyz.openbmc_project.Software.BMC.Updater.service
[Unit]
Description=OpenBMC Software Update Manager

[Service]
ExecStart={sbindir}/phosphor-image-updater
Restart=always
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
