# meta-facebook/meta-santabarbara/recipes-phosphor/gpio/phosphor-gpio-monitor/assert-rmc-main-power-enable.service
[Unit]
Description=RMC Main Power Enabled

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/assert-rmc-main-power-enable

SyslogIdentifier=rmc-main-power-enabled
