# meta-facebook/meta-santabarbara/recipes-phosphor/gpio/phosphor-gpio-monitor/assert-module-power-good@.service
[Unit]
Description=Assert ASIC Module%i Power Good Status

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/assert-module-power-good "%i"

SyslogIdentifier=assert-module%i-power-good
