# meta-facebook/meta-santabarbara/recipes-phosphor/gpio/phosphor-gpio-monitor/deassert-module-power-good@.service
[Unit]
Description=Deassert ASIC Module%i Power Good Status

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/deassert-module-power-good "%i"
SyslogIdentifier=deassert-module%i-power-good