# meta-facebook/meta-santabarbara/recipes-phosphor/gpio/phosphor-gpio-monitor/get-mmc-lastest-version.service
[Unit]
Description=Get the latest MMC version

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/get-mmc-lastest-version

SyslogIdentifier=get-mmc-lastest-version
