# meta-facebook/meta-santabarbara/recipes-phosphor/gpio/phosphor-gpio-monitor/mb-pwr-fault-assert.service
[Unit]
Description=MB Power Fault Assert Handler

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/mb-pwr-fault-handler "assert"
SyslogIdentifier=mb-pwr-fault-assert
