# meta-facebook/meta-santabarbara/recipes-phosphor/gpio/phosphor-gpio-monitor/mctp_remove.service
[Unit]
Description=Remove MCTP Endpoint ID (EID)

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/mctp_remove
