# meta-facebook/meta-santabarbara/recipes-phosphor/gpio/phosphor-gpio-monitor/platform-host-ready.target
[Unit]
Description=santabarbara host ready target
