# meta-facebook/meta-santabarbara/recipes-phosphor/gpio/phosphor-gpio-monitor/swb-pwr-fault-assert.service
[Unit]
Description=SWB Power Fault Assert Handler

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/swb-pwr-fault-handler "assert"
SyslogIdentifier=swb-pwr-fault-assert
