# meta-facebook/meta-santabarbara/recipes-phosphor/gpio/phosphor-gpio-monitor/swb-pwr-fault-deassert.service
[Unit]
Description=SWB Power Fault Deassert Handler

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/swb-pwr-fault-handler "deassert"
SyslogIdentifier=swb-pwr-fault-deassert
