# meta-phosphor/recipes-phosphor/logging/phosphor-logging/xyz.openbmc_project.Syslog.Config.service
[Unit]
Description=Rsyslog config updater

[Service]
ExecStart=/usr/bin/env phosphor-rsyslog-conf
SyslogIdentifier=phosphor-rsyslog-conf
Restart=always
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy=multi-user.target
