# meta-phosphor/recipes-phosphor/logging/phosphor-logging/xyz.openbmc_project.Logging.service
[Unit]
Description=Phosphor Log Manager

[Service]
ExecStart=/usr/bin/env phosphor-log-manager
SyslogIdentifier=phosphor-log-manager
Restart=always
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy=multi-user.target
