# meta-mellanox/meta-msn/recipes-connectivity/mlx-ipmid/mlx-ipmid/mlx_ipmid.service
[Unit]
Description=mlx_ipmid

[Service]
Restart=always
OOMScoreAdjust=-1000

ExecStart=/usr/bin/mlx_ipmid -n
TimeoutSec=120

[Install]
WantedBy=multi-user.target
