# meta-facebook/meta-catalina/recipes-catalina/plat-svc/files/iob-nic-temp-read.service
[Unit]
Description=iob-nic-temp-read
After=xyz.openbmc_project.EntityManager.service
After=xyz.openbmc_project.externalsensor.service
After=mctpd.service

[Service]
Type=simple
ExecStart=/usr/libexec/catalina/iob-nic-temp-read
SyslogIdentifier=iob-nic-temp-read

[Install]
WantedBy=multi-user.target
