# meta-facebook/meta-catalina/recipes-catalina/plat-svc/files/osfp-eeprom-driver-bind.service
[Unit]
Description=Bind OSFP FRU EEPROM driver
After=platform-sys-init.service
After=xyz.openbmc_project.FruDevice.service

[Service]
Type=simple
ExecStart=/usr/libexec/plat-svc/osfp-eeprom-driver-bind
SyslogIdentifier=osfp-eeprom-driver-bind

[Install]
WantedBy=multi-user.target
