# meta-quanta/meta-gsj/recipes-phosphor/fans/phosphor-pid-control/fan-boot-control.service
[Unit]
Description=Set Fan to Default Duty as Booting Up
DefaultDependencies=no

[Service]
Type=simple
ExecStart=/usr/bin/fan-default-speed.sh

[Install]
WantedBy=sysinit.target
