# meta-phosphor/recipes-phosphor/system/obmc-mgr-system/org.openbmc.managers.System.service
[Unit]
Description=Phosphor System Manager
Wants=obmc-mapper.target
After=obmc-mapper.target

[Service]
Restart=always
ExecStart=/usr/bin/env system_manager.py
SyslogIdentifier=system_manager.py
Type=dbus
BusName={BUSNAME}
Environment="PYTHONUNBUFFERED=1"
TimeoutStartSec=180s

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
