# meta-phosphor/common/recipes-phosphor/system/obmc-mgr-system/startup-hacks.service
[Unit]
Description=Startup Hacks

[Service]
Restart=no
Type=oneshot
ExecStart={sbindir}/startup_hacks.sh

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
