# meta-evb/meta-evb-nuvoton/meta-scm-npcm845/recipes-x86/chassis/x86-power-control/obmc-chassis-poweron.target
[Unit]
Description=Chassis Power On
Conflicts=obmc-chassis-poweroff.target