# meta-facebook/meta-yosemite5/meta-rainiera6/recipes-phosphor/mctp/files/rainier-mctp-i3c@.service
[Unit]
Description=Rainier MCTP over I3C setup
Conflicts=host-force-poweroff@%i.service host-graceful-poweroff@%i.service host-powerreset@%i.service obmc-host-stop@%i.target
After=mctpd.service
Before=pldmd.service
Wants=mctpd.service
Requires=mctpd.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/mctp/rainier-mctp-i3c.sh start
ExecStop=/usr/libexec/mctp/rainier-mctp-i3c.sh stop
RemainAfterExit=yes

[Install]
WantedBy=obmc-host-started@%i.target
