# meta-ampere/meta-mitchell/recipes-ampere/platform/ampere-platform-init/ampere-platform-init.service
[Unit]
Description = Ampere Platform Initialization
Before=org.openbmc.control.Power@0.service
Before=xyz.openbmc_project.EntityManager.service

[Service]
Restart=no
RemainAfterExit=true
Type=oneshot
ExecStart=/usr/sbin/ampere_platform_init.sh

[Install]
WantedBy=multi-user.target
