# meta-ampere/meta-mitchell/recipes-ampere/platform/ampere-platform-init/ampere-bmc-heartbeat.service
[Unit]
Description=Ampere BMC heartbeat service

[Service]
Type=simple
Restart=always
ExecStart=/usr/bin/env ampere_bmc_heartbeat.sh
SyslogIdentifier = "ampere-bmc-heartbeat"

[Install]
WantedBy=multi-user.target

