# meta-facebook/meta-yosemite4/recipes-phosphor/gpio/phosphor-gpio-monitor/check-interrupt@.service
[Unit]
Description=Check %i interrupt

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/check-interrupt %i
SyslogIdentifier=check-interrupt%i
