# meta-facebook/meta-yosemite4/recipes-phosphor/gpio/phosphor-gpio-monitor/en-i3c-hub-scan-fru@.service
[Unit]
Description=slot-en-i3c-hub-scan-fru:%i

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/en-i3c-hub-scan-fru %i
SyslogIdentifier=slot-en-i3c-hub-scan-fru %i
