# meta-facebook/meta-yosemite4/recipes-phosphor/gpio/phosphor-gpio-monitor/enable-i3c-hub@.service
[Unit]
Description=slot-enable-i3c-hub:%i

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/enable-i3c-hub %i
SyslogIdentifier=slot-enable-i3c-hub %i
