# meta-facebook/meta-yosemite4/recipes-phosphor/gpio/phosphor-gpio-monitor/fan-board-efuse-fault@.service
[Unit]
Description=fan board power-fault:%i

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/fan-board-efuse-fault %i
SyslogIdentifier=fan-board-efuse-fault%i
