# meta-facebook/meta-yosemite4/recipes-phosphor/gpio/phosphor-gpio-monitor/medusa-hsc-alert@.service
[Unit]
Description=Detect medusa hsc alert on position %i

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/medusa-hsc-alert %i
SyslogIdentifier=medusa-hsc-alert%i
