# meta-facebook/meta-yosemite4/recipes-phosphor/gpio/phosphor-gpio-monitor/medusa-power-module-fault.service
[Unit]
Description=Detect medusa power module fault

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/medusa-power-module-fault
SyslogIdentifier=medusa-power-module-fault
