# meta-facebook/meta-yosemite4/recipes-phosphor/gpio/phosphor-gpio-monitor/nic-power-fault@.service
[Unit]
Description=nic-power-fault

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/nic-power-fault
SyslogIdentifier=nic-power-fault
