# meta-facebook/meta-yosemite4/recipes-phosphor/gpio/phosphor-gpio-monitor/probe-slot-device@.service
[Unit]
Description=Probe Device For Slot%i

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/probe-slot-device %i
SyslogIdentifier=probe-slot-device%i

