# meta-facebook/meta-yosemite4/recipes-phosphor/gpio/phosphor-gpio-monitor/remove-nic-endpoint-slot@.service
[Unit]
Description=Remove MCTP endpoint for Slot%i NIC

[Service]
Type=oneshot
ExecStart=/usr/bin/configure-nic-mctp-endpoint remove %i
SyslogIdentifier=remove-nic-endpoint-slot%i

