# meta-facebook/meta-yosemite4/recipes-phosphor/gpio/phosphor-gpio-monitor/rescan-fru-device@.service
[Unit]
Description=Rescan FruDevice for slot %i
After=probe-slot-device@%i.service
Requires=probe-slot-device@%i.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/rescan-fru-device %i
SyslogIdentifier=rescan-fru-device%i
