# meta-facebook/meta-yosemite4/recipes-phosphor/gpio/phosphor-gpio-monitor/rescan-wf-bic@.service
[Unit]
Description=Rescan WF BIC's EID for slot %i
After=rescan-fru-device@%i.service
Requires=rescan-fru-device@%i.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/rescan-wf-bic %i
SyslogIdentifier=rescan-wf-bic%i
