# meta-facebook/meta-yosemite4/recipes-phosphor/gpio/phosphor-gpio-monitor/setup-nic-endpoint-slot@.service
[Unit]
Description=Setup MCTP endpoint for Slot%i NIC

[Service]
Type=oneshot
ExecStart=/usr/bin/configure-nic-mctp-endpoint add %i
SyslogIdentifier=setup-nic-endpoint-slot%i

