# meta-facebook/meta-yosemite4/recipes-phosphor/gpio/phosphor-gpio-monitor/slot-power-fault@.service
[Unit]
Description=slot power-fault:%i

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/slot-power-fault %i
SyslogIdentifier=slot-power-fault%i
