# meta-facebook/meta-yosemite4/recipes-phosphor/gpio/phosphor-gpio-monitor/slot-hot-plug@.service
[Unit]
Description=slot hot-plug:%i

[Service]
Type=oneshot
ExecStartPre=/usr/libexec/phosphor-state-manager/wait-until-mctp-connection-done %i
ExecStart=/bin/systemctl restart phosphor-discover-system-state@%i.service
