# meta-facebook/meta-yosemite4/recipes-phosphor/gpio/phosphor-gpio-monitor/slot-unplug@.service
[Unit]
Description=slot-unplug %i

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/slot-unplug %i
SyslogIdentifier=slot-unplug %i