# meta-facebook/meta-yosemite4/recipes-phosphor/gpio/phosphor-gpio-monitor/slot-plug-in@.service
[Unit]
Description=slot-plug-in %i

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/slot-plug-in %i
SyslogIdentifier=slot-plug-in %i