# meta-openembedded/meta-webserver/recipes-httpd/sthttpd/sthttpd/thttpd.service
[Unit]
Description=Tiny/Turbo/Throttling Web Server

[Service]
Type=forking
ExecStart=/usr/sbin/thttpd -C @@CONFFILE -c cgi-bin/* -i /var/run/thttpd.pid
PIDFile=/var/run/thttpd.pid

[Install]
WantedBy=multi-user.target
