# import-layers/meta-openembedded/meta-oe/recipes-graphics/xserver-nodm-init/xserver-nodm-init-2.0/xserver-nodm.service
[Unit]
Description=Xserver startup with a display manager

[Service]
EnvironmentFile=/etc/default/xserver-nodm
ExecStart=/etc/X11/Xserver

[Install]
Alias=display-manager.service
