# meta-evb/meta-evb-nuvoton/meta-scm-npcm845/recipes-scm-npcm845/usb-emmc-storage/files/usb_emmc_storage.service
[Unit]
Description=enable usb emmc storage

[Service]
Type=oneshot
ExecStart=/bin/bash -c '/usr/sbin/usb_emmc_storage.sh'

[Install]
WantedBy=multi-user.target
