# meta-openbmc-machines/meta-openpower/meta-ibm/meta-romulus/recipes-phosphor/chassis/avsbus-control/avsbus-disable@.service
[Unit]
Description=Disable the AVS bus on the VRMs
Wants=avsbus-workaround@%i.service
After=avsbus-workaround@%i.service
Before=vrm-control@%i.service
Conflicts=obmc-chassis-poweroff@%i.target
ConditionPathExists=!/run/openbmc/chassis@%i-on

[Service]
ExecStart={bindir}/avsbus-disable.sh
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=obmc-chassis-poweron@%i.target
