# meta-ibm/recipes-phosphor/network/network/ncsi-failover@.service
[Unit]
Description=Stop %i link failover
After=network.target

[Service]
Restart=no
ExecCondition=/lib/systemd/systemd-networkd-wait-online --timeout 30 -i %i
ExecStart=/usr/libexec/ncsi-netlink-ifindex %i --set -p 0 -c 0
SyslogIdentifier=nsci-failover
Type=oneshot

[Install]
# Handled by bitbake recipe
