# meta-facebook/meta-sanmiguel/recipes-phosphor/sensors/dbus-sensors/xyz.openbmc_project.leakdetector.warning.assert@.service
[Unit]
Description=%i Warning Leak Handler

[Service]
Type=oneshot
ExecStart=/usr/libexec/dbus-sensors/warning-leak-handler %i
SyslogIdentifier=%i-warning-leak-handler