# meta-openbmc-machines/meta-openpower/common/recipes-phosphor/fsi/op-fsi/fsi-disable.service
[Unit]
Description=Disable FSI
Conflicts=fsi-enable.service

[Service]
ExecStart=/bin/sh -c "echo fsi-master > /sys/bus/platform/drivers/fsi-master-gpio/unbind"
Type=oneshot
RemainAfterExit=yes
