# meta-facebook/meta-bletchley/recipes-phosphor/state/phosphor-state-manager/chassis-poweron@.service
[Unit]
Description=Turn on AC at Sled%i
Requires=motor-init-calibration@%i.service
After=motor-init-calibration@%i.service
Conflicts=chassis-poweroff@%i.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-state-manager/chassis-poweron %i
SyslogIdentifier=power-ctrl
StartLimitBurst=0

[Install]
WantedBy=obmc-chassis-poweron@%i.target
