# meta-facebook/meta-bletchley/recipes-phosphor/state/phosphor-state-manager/host-powerreset@.service
[Unit]
Description=Sled%i Host Power Reset
Requires=motor-init-calibration@%i.service
After=motor-init-calibration@%i.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/power-ctrl sled%i reset
SyslogIdentifier=power-ctrl
