# meta-phosphor/common/recipes-phosphor/system/obmc-mgr-state/obmc-mgr-state.service
[Unit]
Description=Restore Last State
Wants=obmc-mgr-state.target
After=obmc-mgr-state.target

[Service]
Restart=no
Type=oneshot
ExecStart={sbindir}/discover_system_state.py
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
