# meta-hpe/meta-common/recipes-phosphor/recipes-graphics/obmc-ikvm/obmc-ikvm.service
[Unit]
Description=OpenBMC ipKVM daemon
ConditionPathIsMountPoint=/sys/kernel/config

[Service]
Restart=always
ExecStartPre=/usr/bin/create_usbhid.sh disconnect
ExecStart=/usr/bin/obmc-ikvm -v /dev/video0 -k /dev/hidg0 -p /dev/hidg1 -u 80401000.udc

[Install]
WantedBy=multi-user.target
