# meta-phosphor/recipes-phosphor/interfaces/phosphor-gevent/phosphor-gevent.service
[Unit]
Description=Phosphor Webserver
After=phosphor-gevent.socket
After=obmc-webserver-pre.target

[Service]
Restart=always
ExecStart=/usr/bin/env phosphor-gevent $APPLICATION
SyslogIdentifier=phosphor-gevent
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile={envfiledir}/obmc/wsgi_app
