# meta-ampere/meta-common/recipes-ac01/host/ampere-host-error-monitor/xyz.openbmc_project.ampere_host_error_monitor.service
[Unit]
Description=Ampere Host Error Monitor

[Service]
Restart=always
ExecStart=/usr/bin/ampere-host-error-monitor
Type=simple

[Install]
WantedBy=multi-user.target
