# meta-google/recipes-google/networking/gbmc-bridge/gbmc-br-dhcrelay.service
[Unit]
Description=gBMC DHCP Relay Agent Daemon
After=network.target
StartLimitIntervalSec=10
StartLimitBurst=3

[Service]
Restart=always
RestartSec=5
ExecStart=/usr/libexec/gbmc-br-dhcrelay

[Install]
WantedBy=multi-user.target
