# meta-google/recipes-google/networking/gbmc-bridge/gbmc-br-ra.service
[Unit]
Description=gBMC Bridge RA Discovery
After=network.target
StartLimitIntervalSec=10
StartLimitBurst=3

[Service]
Restart=always
RestartSec=5
ExecStart=/usr/libexec/gbmc-br-ra.sh

[Install]
WantedBy=multi-user.target
