# meta-google/recipes-google/networking/gbmc-bridge/gbmc-br-hostname.service
[Unit]
Before=systemd-networkd.service

[Service]
ExecStart=/usr/libexec/gbmc-br-hostname.sh

[Install]
WantedBy=multi-user.target
