# meta-google/recipes-google/networking/gbmc-bridge/gbmc-br-dhcp@.service
[Unit]
Description=gBMC DHCP Client
After=network.target dhcp-done.service
Wants=dhcp-done.service
StartLimitIntervalSec=10
StartLimitBurst=3

[Service]
Restart=on-failure
RestartSec=5
ExecStartPre=/usr/bin/update-dhcp-status 'ONGOING' 'Starting dhcp process'
ExecStart=/usr/bin/udhcpc6 -f -q -O fqdn -O bootfile_url -O bootfile_param -i %I -s /usr/libexec/gbmc-br-dhcp.sh

[Install]
WantedBy=multi-user.target
