# meta-google/recipes-google/networking/gbmc-bridge/gbmc-br-load-ip.service
[Unit]
Before=gbmc-ip-monitor.service
Before=systemd-networkd.service

[Service]
ExecStart=/bin/bash -c 'source /usr/share/gbmc-br-lib.sh && gbmc_br_reload_ips'

[Install]
WantedBy=multi-user.target
