# meta-google/recipes-google/networking/gbmc-bridge/gbmc-br-ensure-ra.service
[Service]
ExecStart=/usr/libexec/gbmc-br-ensure-ra.sh

[Install]
WantedBy=multi-user.target
