# meta-google/recipes-google/networking/gbmc-bridge/l2-br-dhcp4.service
[Unit]
Description=gBMC L2 Bridge DHCPv4 Client
After=network.target dhcp-done.service
Wants=dhcp-done.service
StartLimitIntervalSec=10
StartLimitBurst=3

[Service]
EnvironmentFile=/usr/share/br-dhcp-env
Restart=always
RestartSec=5
ExecStart=/usr/sbin/udhcpc -f -q -O hostname -O bootfile -i l2br -s /usr/libexec/gbmc-br-dhcp.sh ${L2BR_VENDOR_CLASS_V4}

[Install]
WantedBy=multi-user.target
