# meta-google/recipes-google/host-power-ctrl/gpio-host-pwr/host-ensure-off.service
[Unit]
Description=Ensure host is off on clean AC cycle
Before=host-poweron.service
Before=firmware-updates-pre.target

[Service]
Type=oneshot
ExecStart=/usr/bin/host_ensure_off.sh

[Install]
WantedBy=multi-user.target
