# meta-google/recipes-google/host-power-ctrl/gpio-host-pwr/host-poweron.service
[Unit]
Description=Power on host
After=firmware-updates.target

[Service]
Type=oneshot
ExecStart=/usr/bin/host_poweron.sh

[Install]
WantedBy=multi-user.target
