# meta-google/recipes-google/host-power-ctrl/gpio-host-pwr/host-powercycle.service
[Unit]
Description=Powercycle host

[Service]
Type=oneshot
ExecStart=/usr/bin/host_powercycle.sh
