# meta-google/recipes-google/host-power-ctrl/gpio-host-pwr/host-poweroff-watchdog.service
[Unit]
Description=Power off host

[Service]
Type=oneshot
Environment=DONT_STOP_WATCHDOG=1
ExecStart=/usr/bin/host_poweroff.sh
