# meta-google/recipes-google/host-power-ctrl/gpio-host-pwr/host-reset-cold.service
[Unit]
Description=Cold Reset Host

[Service]
Type=oneshot
ExecStart=/usr/bin/host_reset.sh cold
