# meta-google/recipes-google/host-power-ctrl/gpio-host-pwr/host-poweroff.service
[Unit]
Description=Power off host

[Service]
Type=oneshot
ExecStart=/usr/bin/host_poweroff.sh
