# meta-google/recipes-google/host-power-ctrl/gpio-host-pwr/host-reset-warm.service
[Unit]
Description=Warm Reset Host

[Service]
Type=oneshot
ExecStart=/usr/bin/host_reset.sh warm
