# meta-phosphor/common/recipes-phosphor/ipmi/phosphor-ipmi-host/op-stop-host@.service
[Unit]
Description=Soft power off of the host
Conflicts=obmc-chassis-start@%i.target

[Service]
Restart=no
Type=oneshot
RemainAfterExit=yes
ExecStart={sbindir}/mapper call /org/openbmc/HostServices org.openbmc.HostServices SoftPowerOff

[Install]
WantedBy=obmc-stop-host@%i.target
