# meta-phosphor/common/recipes-phosphor/ipmi/phosphor-ipmi-host/org.openbmc.HostServices.service
[Unit]
Description=Phosphor Inband IPMI
Wants=clear-once.service
After=clear-once.service
Wants=mapper-wait@-org-openbmc-settings-host0.service
After=mapper-wait@-org-openbmc-settings-host0.service

[Service]
Restart=always
ExecStart={sbindir}/ipmid
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
