# meta-phosphor/common/recipes-phosphor/ipmi/phosphor-ipmi-host/phosphor-ipmi-host.service
[Unit]
Description=Phosphor Inband IPMI
# TODO openbmc/openbmc#2059 - The wants/after below should be based on providers
Wants=mapper-wait@-xyz-openbmc_project-control-host0-boot.service
After=mapper-wait@-xyz-openbmc_project-control-host0-boot.service
Wants=mapper-wait@-xyz-openbmc_project-control-host0-boot-one_time.service
After=mapper-wait@-xyz-openbmc_project-control-host0-boot-one_time.service
Wants=mapper-wait@-xyz-openbmc_project-control-host0-power_restore_policy.service
After=mapper-wait@-xyz-openbmc_project-control-host0-power_restore_policy.service
Wants=mapper-wait@-xyz-openbmc_project-control-host0-restriction_mode.service
After=mapper-wait@-xyz-openbmc_project-control-host0-restriction_mode.service
Wants=clear-once.service
After=clear-once.service
# TODO openbmc/openbmc#2008 - remove wants/after on the old settings below
Wants=mapper-wait@-org-openbmc-settings-host0.service
After=mapper-wait@-org-openbmc-settings-host0.service

[Service]
Restart=always
ExecStart=/usr/bin/env ipmid
SyslogIdentifier=ipmid

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
