# meta-openbmc-machines/meta-openpower/common/recipes-phosphor/flash/openpower-software-manager/obmc-flash-bios-updatesymlinks@.service
[Unit]
Description=Updates symlinks for active PNOR version
Wants=mbox-reset.service
Before=mbox-reset.service

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart={sbindir}/obmc-flash-bios updatesymlinks %i
