# meta-openbmc-machines/meta-openpower/common/recipes-phosphor/flash/openpower-software-manager/obmc-flash-bios-ubiremount.service
[Unit]
Description=Remount the squashfs and ubi volumes after a reboot
Wants=org.open_power.Software.Host.Updater.service
Before=org.open_power.Software.Host.Updater.service

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart={sbindir}/obmc-flash-bios ubiremount

[Install]
WantedBy=multi-user.target
