# meta-openbmc-machines/meta-openpower/common/recipes-phosphor/flash/openpower-software-manager/obmc-flash-bios-ubiblock@.service
[Unit]
Description=Create a ubiblock on %I
Requires=obmc-flash-bios-ubiattach.service
After=obmc-flash-bios-ubiattach.service

[Service]
ExecStart=/usr/sbin/ubiblock --create %I

[Install]
WantedBy=multi-user.target
