# meta-openbmc-machines/meta-openpower/common/recipes-phosphor/flash/openpower-software-manager/obmc-flash-bios-ubipatch.service
[Unit]
Description=Create a ubi volume for pnor patches
Wants=usr-local.mount
After=usr-local.mount

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart={sbindir}/obmc-flash-bios ubimount pnor-patch
