# meta-openbmc-machines/meta-openpower/common/recipes-phosphor/flash/openpower-software-manager/obmc-flash-bios-updatesymlinks.service
[Unit]
Description=Updates symlinks for active PNOR version
Before=mboxd.service
Before=mboxd-reload@0.service
After=org.open_power.Software.Host.Updater.service
After=op-reset-chassis-running@0.service
ConditionPathExists=!/run/openbmc/chassis@0-on

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/env obmc-flash-bios updatesymlinks
SyslogIdentifier=obmc-flash-bios

[Install]
WantedBy=mboxd.service mboxd-reload@.service
