# meta-openbmc-machines/meta-openpower/common/recipes-phosphor/flash/openpower-software-manager/obmc-flash-bios-updatesymlinks.service
[Unit]
Description=Updates symlinks for active PNOR version
Before=mboxd.service
Before=mboxd-reload@.service

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart={sbindir}/obmc-flash-bios updatesymlinks

[Install]
WantedBy=mboxd.service mboxd-reload@.service
