# meta-openbmc-machines/meta-openpower/common/recipes-phosphor/flash/openpower-software-manager/obmc-flash-bios-check-clearvolatile@.service
[Unit]
Description=Clear the Volatile PNOR partitions in host%i if Enabled
After=obmc-flash-bios-enable-clearvolatile@%i.service
Before=obmc-host-start-pre@%i.target
ConditionPathExists=!/run/openbmc/host@%i-on

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/env obmc-flash-bios clearvolatile
SyslogIdentifier=obmc-flash-bios-check-clearvolatile

[Install]
RequiredBy=obmc-host-startmin@%i.target
