# meta-openbmc-machines/meta-openpower/common/recipes-phosphor/flash/openpower-software-manager/obmc-flash-bios-ubiclear@.service
[Unit]
Description=Clear UBIFS volume %i

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/env obmc-flash-bios ubiclear %i
SyslogIdentifier=obmc-flash-bios
