# meta-openbmc-machines/meta-openpower/common/recipes-phosphor/flash/openpower-software-manager/obmc-flash-bios-cleanup.service
[Unit]
Description=Clean up UBI volumes unattached to a Host version

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/env obmc-flash-bios ubicleanup
SyslogIdentifier=obmc-flash-bios
