# meta-openbmc-machines/meta-openpower/common/recipes-phosphor/flash/openpower-software-manager/obmc-flash-bios-ubiinit.service
[Unit]
Description=Mount the PNOR UBI volumes

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart={sbindir}/bios-ubiinit
