# meta-openbmc-machines/meta-openpower/common/recipes-phosphor/flash/openpower-software-manager/obmc-flash-bios-ubiattach.service
[Unit]
Description=Attach the PNOR mtd device to UBI

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart={sbindir}/obmc-flash-bios ubiattach
