# meta-openbmc-machines/meta-openpower/common/recipes-phosphor/flash/openpower-software-manager/obmc-flash-bios-ubimount@.service
[Unit]
Description=Mount UBIFS volume pnor-prsv and pnor-rw-%I
Requires=obmc-flash-bios-ubiattach.service
After=obmc-flash-bios-ubiattach.service

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/env obmc-flash-bios ubimount pnor-rw-%i
ExecStart=/usr/bin/env obmc-flash-bios ubimount pnor-prsv
SyslogIdentifier=obmc-flash-bios
