# meta-openbmc-machines/meta-openpower/common/recipes-phosphor/flash/openpower-software-manager/obmc-flash-bios-updatesymlinks.service
[Unit]
Description=Updates symlinks for active PNOR version

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart={sbindir}/obmc-flash-bios updatesymlinks

WantedBy=mboxd.service mboxd-reload@.service
