# meta-bytedance/meta-g220a/recipes-network/network/static-mac-addr/static-mac-addr.service
[Unit]
Description=Enforce Static MAC addr mapping

[Service]
Type=oneshot
Restart=no
ExecStart=/usr/bin/mac-check

[Install]
WantedBy=network.target

