# meta-phosphor/common/recipes-phosphor/host-ipmid/host-ipmid/host-ipmid.service
[Unit]
Description="Phosphor OpenBMC IPMI daemon"

[Service]
ExecStart=/usr/sbin/ipmid

[Install]
WantedBy=multi-user.target
