# meta-phosphor/common/recipes-phosphor/host-ipmid/host-ipmid/host-ipmid.service
[Unit]
Description=Phosphor OpenBMC IPMI daemon
Wants=clear-once.service
After=clear-once.service

[Service]
Restart=always
ExecStart=/usr/sbin/ipmid

[Install]
WantedBy=multi-user.target
