# meta-phosphor/common/recipes-phosphor/host-ipmid/host-ipmid/org.openbmc.HostServices.service
[Unit]
Description=Phosphor Inband IPMI
Wants=clear-once.service
After=clear-once.service

[Service]
Restart=always
ExecStart={sbindir}/ipmid
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
