# meta-facebook/meta-ventura2/recipes-phosphor/state/phosphor-state-manager/chassis-powercycle@.service
[Unit]
Description=Chassis%i Power Cycle

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-state-manager/chassis-powercycle %i

[Install]
WantedBy=obmc-chassis-powercycle@%i.target
