# meta-openbmc-machines/meta-openpower/meta-ibm/meta-witherspoon/recipes-phosphor/chassis/avsbus-control/avsbus-disable@.service
[Unit]
Description=Disable the AVS bus on the VRMs
Wants=power-workarounds@%i.service
After=power-workarounds@%i.service
Before=vrm-control@%i.service
Conflicts=obmc-chassis-poweroff@%i.target
ConditionPathExists=!/run/openbmc/chassis@%i-on

[Service]
ExecStart=/usr/bin/env avsbus-disable.sh
SyslogIdentifier=avsbus-disable.sh
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=obmc-chassis-poweron@%i.target
