# meta-phosphor/recipes-phosphor/chassis/obmc-control-chassis/org.openbmc.control.Chassis@.service
[Unit]
Description=Phosphor Chassis%i Control

[Service]
Restart=always
ExecStart=/usr/bin/env chassis_control.py
SyslogIdentifier=chassis_control.py
Type=dbus
BusName={BUSNAME}
Environment="PYTHONUNBUFFERED=1"
TimeoutStartSec=180s

[Install]
WantedBy=multi-user.target
