# meta-ibm/recipes-phosphor/chassis/vrm-control/ir35221-on-bind@.service
[Unit]
Description=Power on bind ir35221 device driver devs
After=avsbus-enable@%i.service
Wants=obmc-host-start-pre@%i.target
Before=obmc-host-start-pre@%i.target
Conflicts=obmc-chassis-poweroff@%i.target
ConditionPathExists=!/run/openbmc/chassis@%i-on

[Service]
ExecStart=/usr/bin/ir35221-unbind-bind.sh bind
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=obmc-chassis-poweron@%i.target
