# meta-phosphor/common/recipes-phosphor/interfaces/phosphor-gevent/phosphor-gevent-minute-timer.timer
[Unit]
Description=Phosphor Gevent Minute Timer

[Timer]
OnBootSec=1min
OnCalendar=*:0/1
Unit=phosphor-gevent-minute-timer.target

[Install]
WantedBy=basic.target


