# meta-phosphor/common/recipes-phosphor/interfaces/phosphor-gevent/phosphor-gevent.service
[Unit]
Description=Phosphor Webserver
PartOf=systemd-networkd.service
After=phosphor-gevent.socket
After=obmc-webserver-pre.target

[Service]
Restart=always
ExecStart={sbindir}/phosphor-gevent $APPLICATION
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile={envfiledir}/obmc/wsgi_app
