# meta-phosphor/common/recipes-phosphor/interfaces/phosphor-gevent/phosphor-gevent-minute-timer.target
[Unit]
Description=Phosphor Gevent Minute Timer Target
StopWhenUnneeded=yes


