# meta-phosphor/common/recipes-phosphor/interfaces/phosphor-gevent/phosphor-gevent.socket
[Unit]
Description=Phosphor Webserver socket

[Socket]
ListenStream=443

[Install]
WantedBy=sockets.target
