# meta-phosphor/common/recipes-phosphor/host-ipmid/btbridged/btbridged.service
[Unit]
Description=Phosphor OpenBMC BT to DBUS

[Service]
Restart=always
ExecStart=/usr/sbin/btbridged

[Install]
WantedBy=multi-user.target
