# meta-phosphor/common/recipes-phosphor/host-ipmid/btbridged/btbridged.service
[Unit]
Description=Phosphor BT

[Service]
Restart=always
ExecStart={sbindir}/btbridged

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
