# meta-openembedded/meta-multimedia/recipes-mediacenter/kodi/kodi-startup/kodi.service
[Unit]
Description=Kodi media thing

[Service]
User=root
Type=simple
SuccessExitStatus=0 1
IOSchedulingClass=realtime
IOSchedulingPriority=0

ExecStart=/usr/bin/xinit /usr/bin/kodi-standalone -- /usr/bin/X :0 

[Install]
WantedBy=basic.target
