# meta-facebook/meta-yosemite4/recipes-phosphor/state/phosphor-state-manager/chassis-powercycle@.service
[Unit]
Description=Chassis Power Cycle: %i
Conflicts=chassis-poweron@%i.service chassis-poweroff@%i.service host-powercycle@%i.service host-poweroff@%i.service host-poweron@%i.service host-powerreset@%i.service
OnFailure=check-i3c-hub@%i.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-state-manager/chassis-powercycle %i

[Install]
WantedBy=obmc-chassis-powercycle@%i.target
