# meta-facebook/meta-yosemite4/recipes-phosphor/state/phosphor-state-manager/chassis-poweron@.service
[Unit]
Description=power on chassis:%i
Conflicts=chassis-powercycle@%i.service chassis-poweroff@%i.service host-powercycle@%i.service host-poweroff@%i.service host-powerreset@%i.service
OnFailure=chassis-poweron-failure@%i.service check-i3c-hub@%i.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-state-manager/chassis-poweron %i

[Install]
WantedBy=obmc-chassis-poweron@%i.target
