# meta-facebook/meta-yosemite4/recipes-phosphor/state/phosphor-state-manager/chassis-poweron@.service
[Unit]
Description=power on chassis:%i
OnFailure=chassis-poweron-failure@%i.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-state-manager/chassis-poweron %i

[Install]
WantedBy=obmc-chassis-poweron@%i.target
