# meta-facebook/meta-yosemite4/recipes-phosphor/state/phosphor-state-manager/chassis-powercycle@.service
[Unit]
Description=Chassis Power Cycle: %i
Conflicts=host-poweron@%i.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-state-manager/chassis-powercycle %i

[Install]
WantedBy=obmc-chassis-powercycle@%i.target
