# meta-facebook/meta-yosemite4/recipes-phosphor/state/phosphor-state-manager/chassis-poweroff@.service
[Unit]
Description=power off chassis:%i
Conflicts=chassis-powercycle@%i.service chassis-poweron@%i.service host-powercycle@%i.service host-poweroff@%i.service host-poweron@%i.service host-powerreset@%i.service
OnFailure=check-i3c-hub@%i.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-state-manager/chassis-poweroff %i

[Install]
WantedBy=obmc-chassis-hard-poweroff@%i.target
