# meta-facebook/meta-yosemite4/recipes-phosphor/state/phosphor-state-manager/check-i3c-hub@.service
[Unit]
Description=Check chassis%i i3c hub

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-state-manager/check-i3c-hub %i
