# meta-facebook/meta-yosemite4/recipes-phosphor/state/phosphor-state-manager/host-poweron@.service
[Unit]
Description=power on host:%i
Conflicts=chassis-powercycle@%i.service
OnFailure=host-poweron-failure@%i.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-state-manager/host-poweron %i
ExecStartPost=/usr/libexec/phosphor-state-manager/rescan-cxl-eid %i

[Install]
WantedBy=obmc-host-start@%i.target
