# meta-facebook/meta-yosemite4/recipes-phosphor/state/phosphor-state-manager/policy-chassis-poweron@.service
[Unit]
Description=Power Policy for Chassis Power On %i
After=host-poweron@%i.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-state-manager/policy-chassis-poweron %i