# meta-fii/meta-mori/recipes-phosphor/console/obmc-console/host_console_uart_config.service
[Unit]
Description = UART Activation Sequence
Before=host-powerctrl.service

[Service]
ExecStart=/usr/libexec/obmc-console/mori_uart_mux_ctrl.sh
Type=oneshot

[Install]
RequiredBy=multi-user.target
