# meta-intel-openbmc/meta-common/recipes-extended/rsyslog/rsyslog/rotate-event-logs.service
[Unit]
Description=Rotate the event logs

[Service]
Type=simple
ExecStart=/usr/bin/rotate-event-logs.sh

[Install]
WantedBy=multi-user.target
