# meta-ingrasys/meta-zaius/recipes-phosphor/chassis/vcs-control/vcs-off@.service
[Unit]
Description=Turn VCS rails off
Wants=obmc-power-stop-pre@%i.target
Before=obmc-power-stop-pre@%i.target
Conflicts=obmc-host-startmin@%i.target

[Service]
ExecStart=/usr/bin/env zaius_vcs.sh off
SyslogIdentifier=zaius_vcs.sh
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=obmc-host-stop@%i.target
