# upstream-layers/openembedded-core/meta/recipes-connectivity/openssh/openssh/sshd.socket
[Unit]
Conflicts=sshd.service
Wants=sshdgenkeys.service
After=nss-user-lookup.target

[Socket]
ExecStartPre=@BASE_BINDIR@/mkdir -p /var/run/sshd
ListenStream=22
Accept=yes

[Install]
WantedBy=sockets.target
