# meta-phosphor/common/recipes-phosphor/console/obmc-console/obmc-console@.service
[Unit]
Description=Phosphor Console Muxer listening on device /dev/%I
BindsTo=dev-%i.device
After=dev-%i.device

[Service]
# TODO: openbmc/obmc-console#5 - Restore this to 'obmc-console-server ttyVUART0'
ExecStart=/bin/sh -c 'realpath /dev/%i | xargs basename | xargs {sbindir}/obmc-console-server --config {sysconfdir}/obmc-console.conf'
Restart=always

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
