# meta-phosphor/common/recipes-phosphor/console/obmc-console/obmc-console.service
[Unit]
Description=Phosphor Console Muxer

[Service]
ExecStart={sbindir}/obmc-console-server ttyS0
Restart=always

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
