# meta-phosphor/common/recipes-phosphor/console/obmc-console/obmc-console@.service
[Unit]
Description=Phosphor Console Muxer listening on device /dev/%I
BindsTo=dev-%i.device
After=dev-%i.device

[Service]
ExecStart={sbindir}/obmc-console-server --config {sysconfdir}/obmc-console.conf %i
Restart=always

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
