# meta-phosphor/common/recipes-phosphor/console/obmc-console/obmc-console-ssh@.service
[Unit]
Description=Phosphor Host SSH Console
Wants=obmc-console.service

[Service]
Environment="DROPBEAR_RSAKEY_DIR=/etc/dropbear"
EnvironmentFile={envfiledir}/dropbear
ExecStart=-{sbindir}/dropbear -i -r ${DROPBEAR_RSAKEY_DIR}/dropbear_rsa_host_key -c {bindir}/obmc-console-client $DROPBEAR_EXTRA_ARGS
ExecReload={base_bindir}/kill -HUP $MAINPID
StandardInput=socket
KillMode=process
