# meta-phosphor/common/recipes-phosphor/console/obmc-console/obmc-console-ssh@.service
[Unit]
Description=Phosphor Host Console SSH Per-Connection
Wants=obmc-console@{OBMC_CONSOLE_HOST_TTY}.service

[Service]
Environment="DROPBEAR_RSAKEY_DIR=/etc/dropbear"
EnvironmentFile={envfiledir}/dropbear
ExecStart=-/usr/bin/env dropbear -i -r ${{DROPBEAR_RSAKEY_DIR}}/dropbear_rsa_host_key -c "/usr/bin/env obmc-console-client" $DROPBEAR_EXTRA_ARGS
SyslogIdentifier=dropbear
ExecReload={base_bindir}/kill -HUP $MAINPID
StandardInput=socket
KillMode=process
