# meta-phosphor/common/recipes-phosphor/console/obmc-console/obmc-console@.service
[Unit]
Description=Phosphor Console Muxer listening on device /dev/%I
BindsTo=dev-ttyVUART0.device
After=dev-ttyVUART0.device

[Service]
ExecStart={sbindir}/obmc-console-server %i
Restart=always

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
