# meta-openbmc-machines/meta-openpower/meta-ibm/meta-witherspoon/recipes-phosphor/fans/witherspoon-fan-watchdog/witherspoon-fan-watchdog.service
[Unit]
Description=Witherspoon Fan Watchdog Daemon
Wants=phosphor-fan-control-init@%i.service
Before=phosphor-fan-control-init@%i.service
Conflicts=obmc-fan-watchdog-takeover.target

[Service]
EnvironmentFile={envfiledir}/obmc/witherspoon-fan-watchdog/fan-watchdog.conf
ExecStart=/sbin/watchdog -T ${{TIMEOUT}} -t ${{INTERVAL}} -F ${{DEVICE}}
KillSignal=SIGKILL

[Install]
RequiredBy=obmc-chassis-poweron@%i.target
