# upstream-layers/meta-openembedded/meta-networking/recipes-daemons/ippool/ippool/ippool.service
[Unit]
Description=ip address pool allocator
Requires=rpcbind.service
After=rpcbind.service

[Service]
Type=simple
# Start ippoold in the foreground! and enable RPC
ExecStart=@SBINDIR@/ippoold -R -f
# Normal output will go to syslog, so suppress stdout.
StandardOutput=null
# ExecStop is not needed. systemd will send SIGTERM
# and ippoold will exit status 1.
SuccessExitStatus=1

[Install]
WantedBy=multi-user.target

