# meta-evb/meta-evb-nuvoton/meta-scm-npcm845/recipes-scm-npcm845/usb-network/files/usb_network.service
[Unit]
Description=enable usb network
After=phosphor-ipmi-host.service

[Service]
Type=oneshot
ExecStart=/bin/bash -c '/usr/sbin/usb_network.sh'

[Install]
WantedBy=multi-user.target
