# meta-phosphor/recipes-phosphor/mboxd/mboxd/mboxd.service
[Unit]
Description=Phosphor MBOX Daemon
BindsTo=dev-aspeed\x2dlpc\x2dctrl.device
After=dev-aspeed\x2dlpc\x2dctrl.device
Requires=xyz.openbmc_project.Control.Host.service
After=xyz.openbmc_project.Control.Host.service

[Service]
Restart=always
EnvironmentFile=-/etc/default/mboxd
ExecStart=/usr/bin/env mboxd --flash {FLASH_SIZE} --window-size 1M --window-num {WINDOW_NUM} $MBOXD_ARGS
SyslogIdentifier=mboxd

[Install]
WantedBy=multi-user.target
