# meta-phosphor/recipes-phosphor/mboxd/mboxd/mboxd.service
[Unit]
Description=Phosphor MBOX Daemon
BindsTo=dev-aspeed\x2dlpc\x2dctrl.device
After=dev-aspeed\x2dlpc\x2dctrl.device
BindsTo=dev-aspeed\x2dmbox.device
After=dev-aspeed\x2dmbox.device
Requires=phosphor-ipmi-host.service
After=phosphor-ipmi-host.service

[Service]
Restart=always
ExecStart=/usr/bin/env mboxd --flash {FLASH_SIZE} --window-size 1M
SyslogIdentifier=mboxd

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
