# meta-facebook/meta-yosemite4/recipes-phosphor/dbus/phosphor-dbus-monitor/disable-fio-led@.service
[Unit]
Description=Disable slot %i FIO LED

[Service]
Type=oneshot
ExecStart=control-fio-led %i false

