# meta-facebook/meta-yosemite4/recipes-phosphor/dbus/phosphor-dbus-monitor/enable-fio-led@.service
[Unit]
Description=Enable slot %i FIO LED

[Service]
Type=oneshot
ExecStart=control-fio-led %i true

