# meta-facebook/meta-yosemite4/recipes-phosphor/dbus/phosphor-dbus-monitor/fan-rotor-missing-mechanism.service
[Unit]
Description=Service to handle fan rotor missing events

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-dbus-monitor/fan-rotor-missing-mechanism