# meta-facebook/meta-yosemite4/recipes-phosphor/dbus/phosphor-dbus-monitor/fan-rotor-fail-mechanism.service
[Unit]
Description=Add an event log for the fan fail mechanism then shut down 8 hosts.

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-dbus-monitor/fan-rotor-fail-mechanism