# meta-facebook/meta-yosemite4/recipes-phosphor/dbus/phosphor-dbus-monitor/switch-nic-to-main-mode@.service
[Unit]
Description=Switch nic %i to main mode

[Service]
Type=oneshot
ExecStart=switch-nic-main-or-aux-mode %i main

