# meta-phosphor/common/recipes-phosphor/phosphor-time-manager/phosphor-time-manager/org.openbmc.TimeManager.service
[Unit]
Description=Phosphor Time Manager daemon
Wants=mapper-wait@-org-openbmc-NetworkManager-Interface.service
After=mapper-wait@-org-openbmc-NetworkManager-Interface.service
Wants=mapper-wait@-org-openbmc-settings-host0.service
After=mapper-wait@-org-openbmc-settings-host0.service

[Service]
Restart=always
ExecStart={sbindir}/timemanager
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
