# meta-phosphor/common/recipes-phosphor/phosphor-time-manager/phosphor-time-manager/org.openbmc.TimeManager.service
[Unit]
Description=Phosphor Time Manager daemon
Requires=network.service
After=network.service
Requires=settings.service
After=settings.service

[Service]
Restart=always
ExecStart={sbindir}/timemanager
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
