# meta-phosphor/common/recipes-phosphor/system/files/obmc-phosphor-sysd.service
[Unit]
Description=Phosphor OpenBMC system management daemon

[Service]
Restart=always
ExecStart={sbindir}/obmc-phosphor-sysd

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
