# meta-ampere/meta-jade/recipes-phosphor/gpio/ampere-gpio-handling/psu_hotswap_reset@.service
[Unit]
Description=Ampere Computing LLC PSU1 reset hotswap

[Service]
Restart=no
ExecStart=/usr/sbin/ampere_psu_reset_hotswap.sh %i
