# meta-google/recipes-google/ipmi/bare-metal-active/gbmc-bare-metal-active@.target
[Unit]
Description=Untrusted customer OS is running on host CPU %i
