# meta-google/recipes-google/ipmi/bare-metal-active/gbmc-bare-metal-prep@.target
[Unit]
Description=Prepare to start customer OS host CPU %i
