# meta-google/recipes-phosphor/flash/inplace-gbmc-update/inplace-gbmc-verify.service
[Unit]
Description=Verify the Flash Image File

[Service]
Type=oneshot
ExecStart=/usr/bin/inplace-gbmc-verify.sh
