# meta-google/recipes-phosphor/flash/inplace-gbmc-update/inplace-gbmc-version.service
[Unit]
Description=Version string for inplace BMC

[Service]
Type=oneshot
StandardOutput=file:/run/inplace-gbmc-version
StandardError=journal
ExecStartPre=/bin/rm -f /run/inplace-gbmc-version
ExecStart=/usr/bin/inplace-gbmc-version.sh
