# meta-axiado-sdk/recipes-axiado/net/ax-net/eip_fwl.service
[Unit]
Description=eip firewall service
Before=network.target

[Service]
Type=oneshot
ExecStart=/etc/eip/scripts/config_eip_ns.sh

[Install]
WantedBy=multi-user.target
