# poky/meta/recipes-connectivity/nfs-utils/nfs-utils/proc-fs-nfsd.mount
[Unit]
Description=NFSD configuration filesystem
After=systemd-modules-load.service

[Mount]
What=nfsd
Where=/proc/fs/nfsd
Type=nfsd
