# meta-phosphor/common/recipes-phosphor/gpio/phosphor-gpio-monitor/phosphor-gpio-monitor.service
[Unit]
Description=Phosphor GPIO monitor

[Service]
Restart=no
ExecStart={sbindir}/phosphor-gpio-monitor
Type=oneshot
