# meta-quanta/meta-gbs/recipes-google/acpi-power-state/acpi-power-state-daemon/gbs-host-set-boot-failsafe@.service
[Unit]
Description=Set gBMC boot time failsafe
Wants=mapper-wait@-xyz-openbmc_project-state-chassis%i.service
After=mapper-wait@-xyz-openbmc_project-state-chassis%i.service
After=acpi-power-state.service
Before=gbs-sysinit.service

[Service]
Type=exec
ExecStart=/usr/bin/gbs-set-boot-failsafe.sh
RemainAfterExit=yes

[Install]
WantedBy=gbs-host-ready.target
