# meta-quanta/meta-gbs/recipes-google/acpi-power-state/acpi-power-state-daemon/gbs-check-host-state.service
[Unit]
Description = Check Host State to set fan failsafe speed

[Service]
Type=simple
Restart=on-failure
ExecStart=/usr/bin/gbs-check-host-state.sh

[Install]
WantedBy=multi-user.target