# meta-phosphor/common/recipes-phosphor/leds/phosphor-led-manager/obmc-fru-fault-monitor.service
[Unit]
Description=FRU Fault monitor service

[Service]
Restart=yes
ExecStart={sbindir}/phosphor-fru-fault-monitor
Type=oneshot

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
