# meta-phosphor/common/recipes-phosphor/leds/phosphor-led-manager/xyz.openbmc_project.ledmanager.service
[Unit]
Description=Phosphor LED Group Management Daemon
Before=mapper-wait@-xyz-openbmc_project-ledmanager-groups.service

[Service]
Restart=always
ExecStart={sbindir}/phosphor-ledmanager
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
