# meta-facebook/meta-minerva/recipes-phosphor/gpio/phosphor-gpio-monitor/ac-power-good-deassert@.service
[Unit]
Description=%i AC Power Good Status Deassert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/logging "%i AC power good deassert" "BLADE" "%i"
SyslogIdentifier=%i-AC-power-good-deassert
