# meta-facebook/meta-minerva/recipes-phosphor/gpio/phosphor-gpio-monitor/cable-present-deassert@.service
[Unit]
Description=%i Cable Present Deassert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/logging "%i cable present deassert" "CABLE" "%i"
SyslogIdentifier=%i-cable-present-deassert
