# meta-facebook/meta-minerva/recipes-phosphor/gpio/phosphor-gpio-monitor/cable-present-assert@.service
[Unit]
Description=%i Cable Present Assert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/logging "%i cable present assert" "CABLE" "%i"
SyslogIdentifier=%i-cable-present-assert
