# meta-facebook/meta-minerva/recipes-phosphor/gpio/phosphor-gpio-monitor/cr-toggle-boot-enabled.service
[Unit]
Description=Cooling Rack Toggle Boot Enabled

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/cr-toggle-boot-logger LockOutDisabled
SyslogIdentifier=cr-toggleboot-enabled
