# meta-facebook/meta-minerva/recipes-phosphor/gpio/phosphor-gpio-monitor/cr-toggle-boot-disabled.service
[Unit]
Description=Cooling Rack Toggle Boot Disabled

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/cr-toggle-boot-logger "xyz.openbmc_project.State.LockOut.LockOutDisabled" "Critical"
SyslogIdentifier=cr-toggleboot-disabled
