# meta-facebook/meta-minerva/recipes-phosphor/gpio/phosphor-gpio-monitor/fan-reload.service
[Unit]
Description=Reload Fan PID service when Deassertion

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/fan-reload

SyslogIdentifier=fan-reload
