# meta-facebook/meta-minerva/recipes-phosphor/gpio/phosphor-gpio-monitor/leak-detect-deassert@.service
[Unit]
Description=%i Leak Detect Deassert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/logging "%i leak detect deassert" "BLADE" "%i"
SyslogIdentifier=%i-leak-detect-deassert
