# meta-facebook/meta-minerva/recipes-phosphor/gpio/phosphor-gpio-monitor/power-fail-deassert@.service
[Unit]
Description=%i Power Fail Deassert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/logging "%i power fail deassert" "SOURCE" "%i"
SyslogIdentifier=%i-power-fail-deassert
