# meta-facebook/meta-minerva/recipes-phosphor/gpio/phosphor-gpio-monitor/power-fail-assert@.service
[Unit]
Description=%i Power Fail Assert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/logging "%i power fail assert" "SOURCE" "%i"
SyslogIdentifier=%i-power-fail-assert
