# meta-facebook/meta-minerva/recipes-phosphor/gpio/phosphor-gpio-monitor/power-good-deassert@.service
[Unit]
Description=%i Power Good Status Deassert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/logging "%i power good deassert" "BLADE" "%i"
SyslogIdentifier=%i-power-good-deassert
