# meta-facebook/meta-minerva/recipes-phosphor/gpio/phosphor-gpio-monitor/present-deassert@.service
[Unit]
Description=%i Present Deassert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/logging "%i present deassert" "BLADE" "%i"
SyslogIdentifier=%i-present-deassert
