# meta-facebook/meta-minerva/recipes-phosphor/gpio/phosphor-gpio-monitor/present-assert@.service
[Unit]
Description=%i Present Assert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/logging "%i present assert" "BLADE" "%i"
SyslogIdentifier=%i-present-assert
