# meta-facebook/meta-minerva/recipes-phosphor/gpio/phosphor-gpio-monitor/rescan-fru.service
[Unit]
Description=Rescan FruDevice for other blades to support power cycle and hot-plug

[Service]
Type=oneshot
ExecStart=busctl call xyz.openbmc_project.FruDevice /xyz/openbmc_project/FruDevice xyz.openbmc_project.FruDeviceManager ReScan
SyslogIdentifier=rescan-fru-device
