# meta-facebook/meta-minerva/recipes-phosphor/gpio/phosphor-gpio-monitor/rpu-ready-deassert@.service
[Unit]
Description=%i Deassert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/logging "%i deassert" "AALC" "%i"
SyslogIdentifier=rpu-ready-deassert
