# meta-facebook/meta-minerva/recipes-phosphor/gpio/phosphor-gpio-monitor/rpu-ready-deassert.service
[Unit]
Description=RPU ready Deassert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/logging "RPU ready deassert" "AALC" "RPU"
SyslogIdentifier=rpu-ready-deassert
