# meta-facebook/meta-minerva/recipes-phosphor/gpio/phosphor-gpio-monitor/rpu-ready-assert.service
[Unit]
Description=RPU ready Assert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/logging "RPU ready assert" "AALC" "RPU"
SyslogIdentifier=rpu-ready-assert
