# meta-facebook/meta-minerva/recipes-phosphor/gpio/phosphor-gpio-monitor/sfp-present-deassert.service
[Unit]
Description=MGMT SFP Present Deassert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/sfp-present-check 1
SyslogIdentifier=sfp-present-deassert
