# meta-facebook/meta-minerva/recipes-phosphor/gpio/phosphor-gpio-monitor/sfp-present-assert.service
[Unit]
Description=MGMT SFP Present Assert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/sfp-present-check 0
SyslogIdentifier=sfp-present-assert
