# meta-hpe/meta-dl360-g11/recipes-hpe/power-supplies-monitoring/files/psu1-monitoring.service
[Unit]
Description=DL360 power supply 1 monitoring
After=dbus.service

[Service]
Restart=always
RestartSec=5
ExecStart=startMonitoring1.sh

[Install]
WantedBy=multi-user.target
