# meta-openembedded/meta-networking/recipes-support/netperf/files/netserver.service
[Unit]
Description=Netperf Benchmark Server
After=network.target

[Service]
ExecStart=/usr/sbin/netserver -D

[Install]
WantedBy=multi-user.target
