# meta-facebook/meta-greatlakes/recipes-phosphor/state/phosphor-state-manager/host-powercycle@.service
[Unit]
Description=Host Power Cycle: %i

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-state-manager/host-powercycle %i

[Install]
WantedBy=obmc-host-reboot@%i.target
