# upstream-layers/openembedded-core/meta/recipes-graphics/x11-common/xserver-nodm-init/xserver-nodm.service.in
[Unit]
Description=Xserver startup without a display manager
After=systemd-udev-settle.service
Wants=systemd-udev-settle.service

[Service]
EnvironmentFile=/etc/default/xserver-nodm
User=@USER@
ExecStart=/etc/xserver-nodm/Xserver
Restart=always

[Install]
Alias=display-manager.service
