# meta-phosphor/common/recipes-phosphor/fans/phosphor-fan/phosphor-fan-control@.service
[Unit]
Description=Phosphor Fan Control Daemon
Conflicts=obmc-chassis-powered-off@%i.target

[Service]
Restart=on-failure
ExecStart=/usr/bin/env phosphor-fan-control --control
SyslogIdentifier=phosphor-fan-control

[Install]
RequiredBy=obmc-fan-control-ready@%i.target
