# meta-phosphor/common/recipes-phosphor/fans/phosphor-fan/phosphor-fan-monitor@.service
[Unit]
Description=Phosphor Fan Monitor Daemon
Wants=phosphor-fan-control@%i.service
Before=phosphor-fan-control@%i.service
Conflicts=obmc-chassis-poweroff@%i.target

[Service]
Restart=on-failure
ExecStart={sbindir}/phosphor-fan-monitor
Type=notify
NotifyAccess=all

[Install]
RequiredBy=obmc-fan-control-ready@%i.target
