# meta-phosphor/common/recipes-phosphor/fans/phosphor-fan/phosphor-fan-control@.service
[Unit]
Description=Phosphor Fan Control Daemon
Wants=obmc-power-on@%i.target
After=obmc-power-on@%i.target
Conflicts=obmc-power-chassis-off@%i.target

[Service]
Restart=on-failure
ExecStart={sbindir}/phosphor-fan-control

[Install]
RequiredBy=obmc-power-chassis-on@%i.target
