# meta-openbmc-machines/meta-openpower/meta-ibm/meta-romulus/recipes-phosphor/host/op-host-control/vcs_workaround@.service
[Unit]
Description=Run VCS workaround on host%i
Wants=avsbus-enable@%i.service
After=avsbus-enable@%i.service
Conflicts=obmc-power-chassis-off@%i.target

[Service]
ExecStart={sbindir}/openpower-proc-control vcsWorkaround
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=obmc-power-chassis-on@%i.target
