# meta-phosphor/common/recipes-phosphor/fans/obmc-control-fan/obmc-max-fans.service
[Unit]
Description=Set Fans to Maximum
Wants=mapper-wait@-org-openbmc-control-fans.service
After=mapper-wait@-org-openbmc-control-fans.service
After=obmc-fan-control.target

[Service]
Type=oneshot
Restart=no
ExecStart=/bin/sh -c "busctl call `mapper get-service /org/openbmc/control/fans` /org/openbmc/control/fans org.openbmc.control.Fans setMax"
SyslogIdentifier=obmc-max-fans
