# meta-ampere/meta-common/recipes-phosphor/console/obmc-console/obmc-console@.service
[Unit]
Description=%i Console Server
ConditionPathExists=/etc/obmc-console/server.%i.conf

[Service]
# Instance ID is the VUART basename
ExecStartPre=/usr/sbin/ampere_uart_console_setup.sh %i
ExecStart=/usr/sbin/obmc-console-server --config /etc/obmc-console/server.%i.conf %i
SyslogIdentifier=%i-console-server
Restart=always
