# meta-ampere/meta-common/recipes-phosphor/console/obmc-console/obmc-console-ttyS3-ssh.socket
[Unit]
Description=Phosphor Host Console SSH Per-Connection socket
Wants=obmc-console@ttyS3.service

[Socket]
ListenStream=2203
Accept=yes

[Install]
WantedBy=sockets.target
