# meta-phosphor/recipes-phosphor/dump/phosphor-debug-collector/obmc-dump-monitor.service
[Unit]
Description=Phosphor Dump core monitor.

[Service]
ExecStart=/usr/bin/env phosphor-dump-monitor
SyslogIdentifier=phosphor-dump-monitor
Restart=always

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
