# meta-phosphor/recipes-phosphor/dump/phosphor-debug-collector/ramoops-monitor.service
[Unit]
Description=Ramoops monitor.
ConditionPathExistsGlob=/var/lib/systemd/pstore/dmesg-ramoops-*
After=xyz.openbmc_project.Dump.Manager.service

[Service]
Type=oneshot
ExecStart=/usr/bin/phosphor-ramoops-monitor
SyslogIdentifier=ramoops-monitor

[Install]
WantedBy=multi-user.target

