# meta-phosphor/recipes-phosphor/dump/phosphor-debug-collector/xyz.openbmc_project.Dump.Manager.service
[Unit]
Description=Phosphor Dump Manager

[Service]
ExecStartPre=/bin/sh -c 'mkdir -p {BMC_DUMP_PATH}'
ExecStart=/usr/bin/env phosphor-dump-manager
SyslogIdentifier=phosphor-dump-manager
Restart=always
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy=multi-user.target
