# meta-evb/meta-evb-aspeed/meta-evb-ast2500/recipes-phosphor/amd-powerctrl/files/host-gpio.service
[Unit]
Description=Configure GPIOs for Host Power Control

[Service]
Restart=no
RemainAfterExit=true
Type=oneshot
ExecStart=/usr/sbin/init_once.sh
SyslogIdentifier=init_once.sh

[Install]
WantedBy=basic.target
