# meta-evb/meta-evb-aspeed/meta-evb-ast2500/recipes-phosphor/amd-powerctrl/files/host-poweron.service
[Unit]
Description=Start Host
Requires=host-gpio.service
After=host-gpio.service
Conflicts=obmc-chassis-poweroff@i.target

[Service]
RemainAfterExit=yes
Type=oneshot
ExecStart=/usr/bin/env poweron.sh
SyslogIdentifier=poweron.sh

[Install]
WantedBy=obmc-host-start@%i.target
