# meta-phosphor/common/recipes-phosphor/obmc-phosphor-led/obmc-control-led/obmc-led-stop@.service
[Unit]
Description=Turn Off %i LED

[Service]
Type=oneshot
RemainAfterExit=yes
Restart=no
ExecStart={sbindir}/mapper call /org/openbmc/control/led/%i org.openbmc.Led setOff
