# meta-phosphor/common/recipes-phosphor/sensors/obmc-mgr-sensor/org.openbmc.Sensors.service
[Unit]
Description=Phosphor Sensor Manager
Wants=mapper-wait@-xyz-openbmc_project-control-power_supply_redundancy.service
After=mapper-wait@-xyz-openbmc_project-control-power_supply_redundancy.service

[Service]
Restart=always
ExecStart=/usr/bin/env sensor_manager2.py
SyslogIdentifier=sensor_manager2.py
Type=dbus
BusName={BUSNAME}
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
