# meta-quanta/meta-gbs/recipes-gbs/gbs-bios-update/files/phosphor-ipmi-flash-bios-update.service
[Unit]
Description=Phosphor-ipmi-flash update BIOS service

[Service]
ExecStart=/usr/bin/bios-update.sh
Type=oneshot

[Install]
WantedBy=phosphor-ipmi-flash-bios-update.target
