# meta-quanta/meta-gbs/recipes-gbs/gbs-bios-update/files/phosphor-ipmi-flash-bios-verify.service
[Unit]
Description=Phosphor-ipmi-flash verify BIOS service

[Service]
ExecStart=/usr/bin/bios-verify.sh
Type=oneshot

[Install]
WantedBy=phosphor-ipmi-flash-bios-verify.target
