# meta-quanta/meta-gsj/recipes-gsj/quanta-nvme-powerctrl/files/nvme_powermanager.service
[Unit]
Description=SSD NVME Power Manager
Requires=nvme_gpio.service
After=nvme_gpio.service

[Service]
ExecStart=/usr/bin/nvme_powermanager.sh
Restart=always
StandardOutput=syslog

[Install]
WantedBy=multi-user.target
