# meta-phosphor/common/recipes-phosphor/logging/phosphor-logging/org.openbmc.logging.service
[Unit]
Description=Phosphor Settings Manager

[Service]
ExecStart={sbindir}/log_manager
Restart=always
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
