# meta-facebook/meta-yosemite5/recipes-phosphor/gpio/phosphor-gpio-monitor/assert-power-good-drop.service
[Unit]
Description=Assert Power Good Drop

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/assert-power-good-drop

SyslogIdentifier=assert-power-good-drop
