# meta-facebook/meta-yosemite5/recipes-phosphor/gpio/phosphor-gpio-monitor/deassert-host-ready.service
[Unit]
Description=Deassert Host Ready

[Service]
Type=oneshot
ExecStart=busctl set-property xyz.openbmc_project.State.Host0 /xyz/openbmc_project/state/host0 xyz.openbmc_project.State.OperatingSystem.Status OperatingSystemState s xyz.openbmc_project.State.OperatingSystem.Status.OSStatus.Inactive

SyslogIdentifier=deassert-host-ready
