# meta-facebook/meta-yosemite5/recipes-phosphor/gpio/phosphor-gpio-monitor/deassert-power-good-drop.service
[Unit]
Description=Deassert Power Good Drop

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/deassert-power-good-drop

SyslogIdentifier=deassert-power-good-drop
