# meta-facebook/meta-yosemite5/recipes-phosphor/gpio/phosphor-gpio-monitor/fault-led-deassert@.service
[Unit]
Description=%i Power Fault LED Deassert
Conflicts=fault-led-deassert@%i.timer

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/fault-led "%i" "deassert"
SyslogIdentifier=%i-power-fault-led-deassert