# meta-facebook/meta-yosemite5/recipes-phosphor/gpio/phosphor-gpio-monitor/fault-led-assert@.service
[Unit]
Description=%i Power Fault LED Assert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/fault-led "%i" "assert"
SyslogIdentifier=%i-power-fault-led-assert