# meta-facebook/meta-yosemite5/recipes-phosphor/gpio/phosphor-gpio-monitor/gpio_bypass@.service
[Unit]
Description=GPIO BYPASS:%i

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/gpio_bypass %i
SyslogIdentifier=gpio_bypass%i
