# meta-facebook/meta-yosemite5/recipes-phosphor/gpio/phosphor-gpio-monitor/hotplug-ac-cycle@.service
[Unit]
Description=Hotplug AC Cycle trigger:%i

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/hotplug-ac-cycle %i
SyslogIdentifier=hotplug-ac-cycle%i
