# meta-facebook/meta-yosemite5/recipes-phosphor/gpio/phosphor-gpio-monitor/reset_btn@.service
[Unit]
Description=Reset button trigger:%i

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/reset_btn %i
SyslogIdentifier=reset_btn%i
