# meta-facebook/meta-yosemite5/recipes-phosphor/gpio/phosphor-gpio-monitor/thermal-deassert-log@.service
[Unit]
Description=%i Thermal Deassert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/thermal-event-logger "deassert" "%i"
SyslogIdentifier=%i-thermal-deassert
