# meta-facebook/meta-yosemite5/recipes-phosphor/gpio/phosphor-gpio-monitor/thermal-assert-log@.service
[Unit]
Description=%i Thermal Assert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/thermal-event-logger "assert" "%i"
SyslogIdentifier=%i-thermal-assert
