# meta-facebook/meta-yosemite5/recipes-phosphor/gpio/phosphor-gpio-monitor/uart-select-led@.service
[Unit]
Description=Uart Select Led:%i

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/set-uart-select-led %i

SyslogIdentifier=uart-select-led%i
