# yocto-poky/meta/recipes-connectivity/dhcp/files/dhcrelay.service
[Unit]
Description=DHCP Relay Agent Daemon
After=network.target

[Service]
ExecStart=@SBINDIR@/dhcrelay -d --no-pid

[Install]
WantedBy=multi-user.target
