# meta-phosphor/common/recipes-phosphor/chassis/obmc-fsi-scan/fsi-scan.service
[Unit]
Description=Scan FSI devices
Wants=fsi-bind.service
After=fsi-bind.service
Conflicts=obmc-power-chassis-off@%i.target
Conflicts=obmc-fsi-disable.target

[Service]
ExecStart=/bin/sh -c "echo 1 > /sys/bus/platform/devices/fsi-master/scan"
Type=oneshot
RemainAfterExit=yes

[Install]
RequiredBy=obmc-fsi-enable.target
