# meta-fii/meta-kudo/recipes-kudo/network/usb-network/usb_network.service
[Unit]
Description=enable usb network
After=phosphor-ipmi-host.service

[Service]
Type=oneshot
ExecStart=/bin/bash -c '/usr/libexec/usb-network/usb_network.sh'

[Install]
WantedBy=multi-user.target
