# meta-openembedded/meta-oe/recipes-extended/smartmontools/files/smartd.service
[Unit]
Description=Self Monitoring and Reporting Technology (SMART) Daemon
After=syslog.target

[Service]
EnvironmentFile=-@SYSCONFDIR@/smartmontools
ExecStart=@SBINDIR@/smartd -n $smart_opts
ExecReload=@BASE_BINDIR@/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
