# meta-ampere/meta-common/recipes-ampere/host/ac01-boot-progress/ampere-boot-progress.service
[Unit]
Description=Ampere Altra Boot Progress Handling
After=obmc-host-already-on@0.target
BindTo=obmc-host-already-on@0.target

[Service]
ExecStart=/usr/sbin/ampere_boot_progress.sh
SyslogIdentifier=ampere-boot-progress
