# meta-ibm/meta-system1/recipes-phosphor/flash/bios-version/pch-standby.service
[Unit]
Description=Wait for PCH Standby Power Sequence Completion

[Service]
Type=oneshot
ExecStart=/usr/libexec/pch-standby-check.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
