# meta-phosphor/common/recipes-phosphor/obmc-phosphor-event/files/obmc-phosphor-event.service
[Unit]
Description=Phosphor OpenBMC event management daemon

[Service]
Restart=always
ExecStart=/usr/sbin/obmc-phosphor-eventd

[Install]
WantedBy=multi-user.target
