# meta-phosphor/common/recipes-phosphor/obmc-phosphor-event/files/obmc-phosphor-event.service
[Unit]
Description=Phosphor OpenBMC event management daemon

[Service]
EnvironmentFile=-/etc/default/eventd/eventd.conf
Restart=always
ExecStart=/usr/sbin/obmc-phosphor-eventd $OPTIONS

[Install]
WantedBy=multi-user.target
