# meta-phosphor/common/recipes-phosphor/obmc-phosphor-sensor/files/obmc-phosphor-sensord.service
[Unit]
Description=Phosphor OpenBMC sensor management daemon

[Service]
ExecStart=/usr/sbin/obmc-phosphor-sensord

[Install]
WantedBy=multi-user.target
