# meta-facebook/meta-harma/recipes-phosphor/gpio/phosphor-gpio-monitor/assert-power-good.service
[Unit]
Description=Assert Power Good Status

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/assert-power-good

SyslogIdentifier=assert-power-good
