# meta-facebook/meta-harma/recipes-phosphor/gpio/phosphor-gpio-monitor/assert-gpio-log@.service
[Unit]
Description=%i GPIO Assert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/logging "%i gpio assert" "GPIO" "%i"
SyslogIdentifier=%i-gpio-assert
