# meta-facebook/meta-harma/recipes-phosphor/gpio/phosphor-gpio-monitor/auto-poweroff@.service
[Unit]
Description=auto power off

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/auto-poweroff %i

SyslogIdentifier=auto-poweroff
