# meta-facebook/meta-harma/recipes-phosphor/gpio/phosphor-gpio-monitor/auto-poweron@.service
[Unit]
Description=auto power on
After=device-reinitial@.service
Wants=device-reinitial@.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/auto-poweron %i

SyslogIdentifier=auto-power-on
