# meta-facebook/meta-harma/recipes-phosphor/gpio/phosphor-gpio-monitor/deassert-power-good.service
[Unit]
Description=Deassert Power Good Status

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/deassert-power-good

SyslogIdentifier=deassert-power-good
