# meta-facebook/meta-harma/recipes-phosphor/gpio/phosphor-gpio-monitor/deassert-uart-switch-button.service
[Unit]
Description=Deassert Uart Switch Button

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/deassert-uart-switch-button

SyslogIdentifier=deassert-uart-switch-button
