# meta-facebook/meta-harma/recipes-phosphor/gpio/phosphor-gpio-monitor/deassert-gpio-log@.service
[Unit]
Description=%i GPIO Deassert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/logging "%i gpio deassert" "GPIO" "%i"
SyslogIdentifier=%i-gpio-deassert
