# meta-facebook/meta-harma/recipes-phosphor/gpio/phosphor-gpio-monitor/deassert-post-end.service
[Unit]
Description=Deassert Post Complete Status

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/deassert-post-end

SyslogIdentifier=deassert-post-end
