# meta-facebook/meta-harma/recipes-phosphor/gpio/phosphor-gpio-monitor/deassert-reset-button.service
[Unit]
Description=Deassert Reset Button

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/deassert-reset-button

SyslogIdentifier=deassert-reset-button
