# meta-facebook/meta-harma/recipes-phosphor/gpio/phosphor-gpio-monitor/device-reinitial.service
[Unit]
Description=Device Reinitial

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/device-reinitial

SyslogIdentifier=device-reinitial
