# meta-facebook/meta-harma/recipes-phosphor/gpio/phosphor-gpio-monitor/gpios-assert-log@.service
[Unit]
Description=%i GPIO Assert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/gpios-event-logger "assert" "%i"
SyslogIdentifier=%i-gpio-assert
