# meta-facebook/meta-harma/recipes-phosphor/gpio/phosphor-gpio-monitor/gpios-deassert-log@.service
[Unit]
Description=%i GPIO Deassert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/gpios-event-logger "deassert" "%i"
SyslogIdentifier=%i-gpio-deassert
