# meta-facebook/meta-harma/recipes-phosphor/gpio/phosphor-gpio-monitor/host-power-on.service
[Unit]
Description=Host power on assertion

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/host-power-on

SyslogIdentifier=host-power-on
