# meta-facebook/meta-harma/recipes-phosphor/gpio/phosphor-gpio-monitor/host-power-off.service
[Unit]
Description=Host power off assertion

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/host-power-off

SyslogIdentifier=host-power-off
