# meta-facebook/meta-harma/recipes-phosphor/gpio/phosphor-gpio-monitor/mmc-recovery.service
[Unit]
Description=mmc recovery

[Service]
Type=oneshot
ExecStart=/usr/libexec/mctp/setup-devices-eid

SyslogIdentifier=mmc-recovery
