# meta-facebook/meta-harma/recipes-phosphor/gpio/phosphor-gpio-monitor/mmc-recovery.service
[Unit]
Description=mmc recovery service

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/mmc-recovery

SyslogIdentifier=mmc-recovery
