# meta-facebook/meta-harma/recipes-phosphor/gpio/phosphor-gpio-monitor/multi-gpios-sys-init.service
[Unit]
Wants=mapper-wait@-xyz-openbmc_project-state-host0.service
After=mapper-wait@-xyz-openbmc_project-state-host0.service


[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/multi-gpios-sys-init
SyslogIdentifier=Power Good Status Init

[Install]
WantedBy=multi-user.target
