# meta-facebook/meta-harma/recipes-phosphor/gpio/phosphor-gpio-monitor/smc-assert-log@.timer
[Unit]
Description=Timer to trigger SMC ready
[Timer]
OnActiveSec=15s
Unit=smc-assert-log@%i.service
[Install]
WantedBy=timers.target
