# meta-facebook/meta-harma/recipes-phosphor/gpio/phosphor-gpio-monitor/smc-assert-log@.service
[Unit]
Description=%i SMC Assert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/smc-event-logger "assert" "%i"
SyslogIdentifier=%i-smc-assert
