# meta-facebook/meta-harma/recipes-phosphor/gpio/phosphor-gpio-monitor/smc-deassert-log@.service
[Unit]
Description=%i SMC Deassert

[Service]
Type=oneshot
ExecStartPre=/bin/systemctl stop smc-assert-log@%i.timer
ExecStart=/usr/libexec/phosphor-gpio-monitor/smc-event-logger "deassert" "%i"
SyslogIdentifier=%i-smc-deassert

Conflicts=smc-assert-log@%i.timer
