# meta-facebook/meta-harma/recipes-phosphor/gpio/phosphor-gpio-monitor/thermaltrip-deassert-log@.service
[Unit]
Description=%i Thermal Trip Deassert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/thermal-event-logger DeviceOperatingNormalTemperature %i
SyslogIdentifier=%i-thermaltrip-deassert
