# meta-facebook/meta-harma/recipes-phosphor/gpio/phosphor-gpio-monitor/thermtrip-assert-log.service
[Unit]
Description=Thermtrip Assert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/thermal-event-logger assert cpu0-thermtrip-alert
SyslogIdentifier=thermtrip-assert
