# meta-facebook/meta-harma/recipes-phosphor/gpio/phosphor-gpio-monitor/thermtrip-deassert-log.service
[Unit]
Description=Thermtrip Deassert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/thermal-event-logger deassert cpu0-thermtrip-alert
SyslogIdentifier=thermtrip-deassert
