# import-layers/meta-openembedded/meta-networking/recipes-daemons/proftpd/files/proftpd.service
[Unit]
Description=proftpd Daemon

[Service]
Type=forking
ExecStart=-@SBINDIR@/proftpd -c @SYSCONFDIR@/proftpd.conf
StandardError=syslog
