# meta-google/recipes-google/nftables/files/nftables.service
[Unit]
Before=network-pre.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/nft-configure.sh
ExecStop=/bin/bash -c 'nft flush ruleset'

[Install]
WantedBy=multi-user.target
