# common/recipes-phosphor/rest-dbus/files/rest-dbus.service
[Unit]
Description=DBUS introspecting REST server.

[Service]
ExecStart=/usr/bin/rest-dbus

[Install]
WantedBy=multi-user.target
