# meta-ibm/meta-witherspoon/recipes-phosphor/chassis/avsbus-control/swift/avsbus-enable@.service
[Unit]
Description=Enable the AVS bus on VRMs
Wants=op-power-start@%i.service
Before=op-power-start@%i.service
Conflicts=obmc-chassis-poweroff@%i.target
ConditionPathExists=!/run/openbmc/chassis@%i-on

[Service]
ExecStart={bindir}/avsbus-enable.sh
SyslogIdentifier=avsbus-enable.sh
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=obmc-chassis-poweron@%i.target
