# poky/meta/files/overlayfs-create-dirs.service.in
[Unit]
Description=Overlayfs directories setup {LOWERDIR}
RequiresMountsFor={DATA_MOUNT_POINT}
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/sbin/overlayfs-create-dirs.sh {LOWERDIR} {DATA_MOUNT_POINT}
RemainAfterExit=true
StandardOutput=journal

[Install]
WantedBy=multi-user.target
