# poky/meta/files/overlayfs-create-dirs.service.in
[Unit]
Description=Overlayfs directories setup
Requires={DATA_MOUNT_UNIT}
After={DATA_MOUNT_UNIT}
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=mkdir -p {DATA_MOUNT_POINT}/upper{LOWERDIR}
ExecStart=mkdir -p {DATA_MOUNT_POINT}/workdir{LOWERDIR}
RemainAfterExit=true
StandardOutput=journal

[Install]
WantedBy=multi-user.target
