# meta-phosphor/recipes-phosphor/leds/phosphor-led-manager/obmc-fru-fault-monitor.service
[Unit]
Description=FRU Fault monitor service
Wants=mapper-wait@-xyz-openbmc_project-led-groups.service
After=mapper-wait@-xyz-openbmc_project-led-groups.service

[Service]
Restart=always
ExecStart=/usr/bin/env phosphor-fru-fault-monitor
SyslogIdentifier=phosphor-fru-fault-monitor

[Install]
WantedBy=multi-user.target
