# meta-phosphor/recipes-phosphor/leds/phosphor-led-manager/xyz.openbmc_project.LED.GroupManager.service
[Unit]
Description=Phosphor LED Group Management Daemon
Before=mapper-wait@-xyz-openbmc_project-led-groups.service
Wants=obmc-mapper.target
After=obmc-mapper.target

[Service]
Restart=always
ExecStart=/usr/bin/env phosphor-ledmanager
SyslogIdentifier=phosphor-ledmanager
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy=multi-user.target
