# meta-phosphor/common/recipes-phosphor/obmc-console/obmc-console/obmc-console.service
[Unit]
Description=OpenBMC console daemon

[Service]
ExecStart=/usr/sbin/obmc-console-server ttyS0
Restart=always

[Install]
WantedBy=multi-user.target
