# meta-phosphor/recipes-phosphor/leds/phosphor-led-sysfs/xyz.openbmc_project.led.controller@.service
[Unit]
Description=Phosphor sysfs LED controller

[Service]
Restart=always
ExecStart=/usr/bin/env phosphor-ledcontroller -p %f
SyslogIdentifier=phosphor-ledcontroller
