# import-layers/meta-openembedded/meta-efl/recipes-efl/efl/entrance/entrance.service
[Unit]
Description=Login manager for Enlightenment

[Service]
ExecStart=/usr/sbin/entrance

[Install]
Alias=display-manager.service
