# meta-google/recipes-phosphor/console/obmc-console/host-console-journal.service.in
[Unit]
Requires=obmc-console@@HOST_TTY@.service
After=obmc-console@@HOST_TTY@.service

[Service]
Type=simple
SyslogIdentifier=host-console
ExecStart=/usr/bin/host-console-filtered.sh

[Install]
WantedBy=multi-user.target
