# meta-google/recipes-phosphor/console/obmc-console/serial-to-host@.service.in
[Unit]
Description=Host Serial Getty on %I
BindsTo=dev-%i.device
After=dev-%i.device systemd-user-sessions.service
After=rc-local.service
Before=getty.target
IgnoreOnIsolate=yes

Requires=obmc-console@@HOST_TTY@.service
After=obmc-console@@HOST_TTY@.service
Conflicts=serial-to-bmc@%i.service
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=0
RestartPreventExitStatus=2
SuccessExitStatus=2
ExecStart=-/sbin/agetty -8 -n -l /usr/bin/obmc-console-client -o '-e a00pae9n2B9m1c' -i -N -R -L %I 115200 xterm
ExecStopPost=-/bin/sh -c '(( EXIT_STATUS == 2 )) && systemctl start --no-block serial-to-bmc@%i'
UtmpIdentifier=%I
TTYPath=/dev/%I
KillMode=process
IgnoreSIGPIPE=no
SendSIGHUP=yes

[Install]
WantedBy=getty.target
