# meta-facebook/meta-bletchley/recipes-bletchley/detect-gpio-present/files/bletchley-sled-removal@.target
[Unit]
Description=Bletchlet SLED%i removal target

Conflicts=motor-init-calibration@%i.service
Conflicts=host-ac-on@%i.service
Conflicts=bletchley-sled-insertion@%i.target

RefuseManualStop=yes