# meta-openbmc-machines/meta-openpower/common/recipes-phosphor/inventory-upload/inventory-upload/op-inventory-upload@.service
[Unit]
Description=Inventory Upload
Wants=obmc-host-start-pre@%i.target
Before=obmc-host-start-pre@%i.target
Conflicts=obmc-host-stop@%i.target

[Service]
Restart=no
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/env inventory_upload.py
SyslogIdentifier=inventory_upload.py
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=obmc-host-start@%i.target
