# meta-facebook/meta-ventura/recipes-phosphor/gpio/phosphor-gpio-monitor/leak-detect-assert@.service
[Unit]
Description=%i Leak Detect Assert

[Service]
Type=oneshot
ExecStart=/usr/bin/log-create xyz.openbmc_project.State.Leak.Detector.LeakDetectedWarning --json "{\"DETECTOR_NAME\": \"/xyz/openbmc_project/state/leak/detector/%i\"}"
SyslogIdentifier=%i-small-leak-detect-assert
