# meta-facebook/meta-ventura/recipes-phosphor/gpio/phosphor-gpio-monitor/led-amber-deassert@.service
[Unit]
Description=%i Leak Detect Deassert
[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/frontled "%i leak detect deassert" "BLADE" "%i"
SyslogIdentifier=%i-leak-detect-deassert
