# meta-facebook/meta-ventura/recipes-phosphor/gpio/phosphor-gpio-monitor/led-blue-init-deassert@.service
[Unit]
Description=%i Power Good LED Init Deassert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/frontled "%i" "true" "deassert"
SyslogIdentifier=%i-power-good-led-init-deassert
