# meta-facebook/meta-ventura/recipes-phosphor/gpio/phosphor-gpio-monitor/led-blue-deassert@.service
[Unit]
Description=%i Power Good LED Deassert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/frontled "%i" "false" "deassert"
SyslogIdentifier=%i-power-good-led-deassert
