# meta-facebook/meta-ventura/recipes-phosphor/gpio/phosphor-gpio-monitor/led-blue-init-assert@.service
[Unit]
Description=%i Power Good LED Init Assert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/frontled "%i" "true" "assert"
SyslogIdentifier=%i-power-good-led-init-assert
