# meta-facebook/meta-ventura/recipes-phosphor/gpio/phosphor-gpio-monitor/led-blue-assert@.service
[Unit]
Description=%i Power Good LED Assert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/frontled "%i" "power good LED" "assert"
SyslogIdentifier=%i-power-good-led-assert
