# meta-facebook/meta-ventura/recipes-phosphor/gpio/phosphor-gpio-monitor/led-amber-assert@.service
[Unit]
Description=%i Leak Detect Assert
[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/frontled "%i leak detect assert" "BLADE" "%i"
SyslogIdentifier=%i-leak-detect-assert
