# meta-facebook/meta-ventura/recipes-phosphor/gpio/phosphor-gpio-monitor/power-good-assert@.service
[Unit]
Description=%i Power Good Status Assert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/logging "%i power good assert" "BLADE" "%i"
SyslogIdentifier=%i-power-good-assert
