# meta-facebook/meta-ventura/recipes-phosphor/gpio/phosphor-gpio-monitor/rack-level-leak@.service
[Unit]
Description=%i Rack Level Leak

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/rack-level-leak
SyslogIdentifier=%i-rack-level-leak
