# meta-facebook/meta-ventura/recipes-phosphor/gpio/phosphor-gpio-monitor/smc-deassert-log@.service
[Unit]
Description=%i SMC Deassert

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/smc-event-logger "deassert" "%i"
SyslogIdentifier=%i-smc-deassert
