# meta-facebook/meta-ventura/recipes-phosphor/gpio/phosphor-gpio-monitor/tray-identify@.service
[Unit]
Description=Tray Identify Service for %i
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/bin/tray-identify %i

[Install]
WantedBy=multi-user.target