# meta-facebook/recipes-fb/hmc-datetime-sync/files/hmc-datetime-sync.service
[Unit]
Description=Sync HMC DateTime via Redfish

[Service]
Type=simple
ExecStart=/usr/libexec/hmc-datetime-sync/hmc-datetime-sync

[Install]
WantedBy=multi-user.target
