# meta-facebook/recipes-fb/hmc-datetime-sync/files/hmc-datetime-sync.timer
[Unit]
Description=Periodic HMC DateTime Sync Timer

[Timer]
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
