# poky/meta/recipes-core/volatile-binds/files/volatile-binds.service.in
[Unit]
Description=Bind mount volatile @where@
DefaultDependencies=no
Before=local-fs.target
RequiresMountsFor=@whatparent@ @whereparent@
ConditionPathIsReadWrite=@whatparent@
ConditionPathExists=@where@
ConditionPathIsReadWrite=!@where@

[Service]
Type=oneshot
RemainAfterExit=Yes
TimeoutSec=0
Environment=MOUNT_COPYBIND_AVOID_OVERLAYFS=@avoid_overlayfs@
ExecStart=/sbin/mount-copybind @what@ @where@
ExecStop=/bin/umount @where@

[Install]
WantedBy=local-fs.target
