# meta-bytedance/meta-g220a/recipes-phosphor/datetime/me-time-sync/me-time-sync.service
[Unit]
Description=Sync time from ME
Wants=ipmb.service
After=ipmb.service

[Service]
Type=oneshot
ExecStart=/usr/bin/me-time-sync.sh

[Install]
WantedBy=multi-user.target
