# meta-ibm/meta-witherspoon/recipes-phosphor/network/network/ncsi-netlink.service
[Unit]
Description=Stop the ethernet link failover
Wants=xyz.openbmc_project.Network.service
After=xyz.openbmc_project.Network.service

[Service]
Restart=no
ExecStart=/usr/bin/env ncsi-netlink --set -x 2 -p 0 -c 0
SyslogIdentifier=ncsi-netlink
Type=oneshot

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
