# meta-axiado-sdk/recipes-axiado/sysproxy/sysproxy/sysproxy.service
[Unit]
Description=Sysproxy Service
After=eip_fwl.service

[Service]
ExecStart=/usr/bin/sysmgr_proxy
Restart=always
StandardOutput=tty
TTYPath=/dev/ttyPS3

[Install]
WantedBy=multi-user.target
