# meta-facebook/meta-anacapa/recipes-phosphor/gpio/phosphor-gpio-monitor/dimm-power-fault-logger@.service
[Unit]
Description=DIMM Power Fault Log

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/dimm-power-fault-logger "%i"
SyslogIdentifier=dimm-power-fault-logger
