# meta-facebook/meta-anacapa/recipes-phosphor/gpio/phosphor-gpio-monitor/fan-pg-event-logger@.service
[Unit]
Description=Fan Power Good Event Logger

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/fan-pg-event-logger "%i"
SyslogIdentifier=fan-pg-event-logger