# meta-facebook/meta-anacapa/recipes-phosphor/gpio/phosphor-gpio-monitor/ssd-pg-event-logger@.service
[Unit]
Description=SSD Power Good Event Logger for %I
After=xyz.openbmc_project.State.Chassis.service
Wants=xyz.openbmc_project.State.Chassis.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/ssd-pg-event-logger "%i"