# meta-ampere/meta-jade/recipes-phosphor/gpio/psu-hotswap-reset/psu1_hotswap_reset.service
[Unit]
Description=Ampere Computing LLC PSU1 reset hotswap

[Service]
Restart=no
ExecStart=ampere_psu_reset_hotswap.sh 1
