# meta-ibm/recipes-phosphor/fans/phosphor-fan/witherspoon/phosphor-fan-control@.service
[Unit]
Description=Phosphor Fan Control Daemon
Requires=max31785-hwmon-helper@ahb-apb-bus\x401e78a000-i2c\x2dbus\x40{ADDR}-max31785\x4052.service
After=max31785-hwmon-helper@ahb-apb-bus\x401e78a000-i2c\x2dbus\x40{ADDR}-max31785\x4052.service

[Service]
Restart=on-failure
ExecStart=/usr/bin/phosphor-fan-control --control

[Install]
RequiredBy=obmc-chassis-poweron@%i.target
