# meta-arm/meta-arm-bsp/recipes-security/optee/optee-client/tee-supplicant@.service
[Unit]
Description=TEE Supplicant on %i

[Service]
User=root
EnvironmentFile=-@sysconfdir@/default/tee-supplicant
ExecStart=@sbindir@/tee-supplicant $OPTARGS

[Install]
WantedBy=basic.target
