# meta-openembedded/meta-oe/recipes-bsp/firmwared/firmwared/firmwared.service
[Unit]
Description=Linux Firmware Loader Daemon

[Service]
Type=simple
ExecStart=@BINDIR@/firmwared -d /lib/firmware

[Install]
WantedBy=multi-user.target
