# meta-quanta/meta-olympus-nuvoton/recipes-x86/chassis/x86-power-control/obmc-chassis-poweroff.target
[Unit]
Description=Chassis Power Off
Conflicts=obmc-chassis-poweron.target
