# meta-nuvoton/recipes-nuvoton/program-edid/program-edid/program-edid.service
[Unit]
Description=Program edid service

[Service]
Restart=no
Type=oneshot
ExecStart=/usr/bin/program-edid.sh ${{SOURCE}} ${{TARGET}}
SyslogIdentifier=program-edid
EnvironmentFile=/etc/default/obmc/edid/program_edid

[Install]
WantedBy=multi-user.target