# yocto-poky/meta/recipes-graphics/wayland/weston-init/weston.service
[Unit]
Description=Weston Wayland compositor startup
RequiresMountsFor=/run

[Service]
User=root
EnvironmentFile=-/etc/default/weston
Environment="XDG_RUNTIME_DIR=/run/user/root"
ExecStartPre=/bin/mkdir -p /run/user/root
ExecStartPre=/bin/chmod 0700 /run/user/root
ExecStart=/usr/bin/openvt -v -e /usr/bin/weston -- $OPTARGS

[Install]
WantedBy=multi-user.target

