# meta-facebook/meta-yosemite5/meta-rainiera6/recipes-phosphor/gpio/phosphor-gpio-monitor/rainier-gpio-handler@.service
[Unit]
Description=Rainier GPIO handler (%i)

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-gpio-monitor/rainier-gpio-handler %i
SyslogIdentifier=rainier-gpio-handler-%i
