# meta-phosphor/common/recipes-phosphor/network/network/network-update-dns.service
[Unit]
Description=Phosphor DNS Settings Cache Manager
Wants=mapper-wait@-org-openbmc-NetworkManager-Interface.service
After=mapper-wait@-org-openbmc-NetworkManager-Interface.service

[Service]
ExecStart={sbindir}/netman_watch_dns /run/systemd/netif/state
Restart=always

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
