# meta-phosphor/common/recipes-phosphor/network/network/network-update-dns.service
[Unit]
Description=Phosphor DNS Updater
Requires=network.service
After=network.service

[Service]
ExecStart={sbindir}/netman_watch_dns /run/systemd/netif/state
Restart=always

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
