# meta-phosphor/common/recipes-phosphor/network/network/network.service
[Unit]
Description=Network DBUS object
Requires=skeleton.service
After=skeleton.service

[Service]
ExecStart=/usr/sbin/netman.py

[Install]
WantedBy=multi-user.target
