# meta-phosphor/common/recipes-phosphor/network/network/network-update-dns.service
[Unit]
Description=DNS update handler
Requires=network.service
After=network.service

[Service]
ExecStart=/usr/sbin/netman_watch_dns /run/systemd/netif/state

[Install]
WantedBy=multi-user.target
