# meta-phosphor/common/recipes-phosphor/network/network/network-update-dns.service
[Unit]
Description=Phosphor DNS Settings Cache Manager
Wants=mapper-wait@-org-openbmc-NetworkManager-Interface.service
After=mapper-wait@-org-openbmc-NetworkManager-Interface.service

[Service]
ExecStart=/usr/bin/env netman_watch_dns /run/systemd/netif/state
SyslogIdentifier=netman_watch_dns
Restart=always

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
