# meta-phosphor/common/recipes-phosphor/network/network/network.service
[Unit]
Description=Network DBUS object

[Service]
ExecStart=/usr/sbin/netman.py
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
