# meta-phosphor/common/recipes-phosphor/network/network/org.openbmc.NetworkManager.service
[Unit]
Description=Phosphor Network Manager

[Service]
ExecStart=/usr/bin/env netman.py
SyslogIdentifier=netman.py
Environment="PYTHONUNBUFFERED=1"
Restart=always
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
