# meta-phosphor/common/recipes-phosphor/network/network/xyz.openbmc_project.Network.service
[Unit]
Description=Phosphor Network Manager
After=obmc-mapper.target

[Service]
ExecStart={sbindir}/phosphor-network-manager
Restart=always
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}

