# meta-facebook/meta-catalina/recipes-phosphor/state/phosphor-state-manager/chassis-powercycle@.service
[Unit]
Description=Chassis Power Cycle: %i
Before=obmc-chassis-powercycle@%i.target

[Service]
Type=oneshot
ExecStart=/usr/libexec/phosphor-state-manager/chassis-powercycle %i

[Install]
RequiredBy=obmc-chassis-powercycle@%i.target
