# import-layers/yocto-poky/meta/recipes-graphics/wayland/weston-init/weston.service
[Unit]
Description=Weston Wayland Compositor
RequiresMountsFor=/run

[Service]
User=root
EnvironmentFile=-/etc/default/weston
ExecStart=/usr/bin/weston-start -v -e -- $OPTARGS

[Install]
WantedBy=multi-user.target

