# meta-intel/meta-common/recipes-extended/rsyslog/rsyslog/rotate-event-logs.timer
[Unit]
Description=Run rotate-event-logs.service every minute

[Timer]
OnActiveSec=60
OnUnitActiveSec=60

[Install]
WantedBy=timers.target
