# meta-phosphor/common/recipes-phosphor/syslogmgmt/obmc-syslog-mgmt/org.openbmc.SyslogManagement.service
[Unit]
Description=Phosphor Syslog settings manager

[Service]
Restart=always
ExecStart={sbindir}/syslogmgmt.py
Type=dbus
BusName={BUSNAME}
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
