# meta-facebook/meta-bletchley/recipes-bletchley/plat-svc/files/bletchley-host-state-monitor.service
[Unit]
Description=Bletchley host state monitoring
Wants=mapper-wait@-xyz-openbmc_project-state-chassis1.service
After=mapper-wait@-xyz-openbmc_project-state-chassis1.service
Wants=mapper-wait@-xyz-openbmc_project-state-chassis2.service
After=mapper-wait@-xyz-openbmc_project-state-chassis2.service
Wants=mapper-wait@-xyz-openbmc_project-state-chassis3.service
After=mapper-wait@-xyz-openbmc_project-state-chassis3.service
Wants=mapper-wait@-xyz-openbmc_project-state-chassis4.service
After=mapper-wait@-xyz-openbmc_project-state-chassis4.service
Wants=mapper-wait@-xyz-openbmc_project-state-chassis5.service
After=mapper-wait@-xyz-openbmc_project-state-chassis5.service
Wants=mapper-wait@-xyz-openbmc_project-state-chassis6.service
After=mapper-wait@-xyz-openbmc_project-state-chassis6.service
Wants=mapper-wait@-xyz-openbmc_project-state-host1.service
After=mapper-wait@-xyz-openbmc_project-state-host1.service
Wants=mapper-wait@-xyz-openbmc_project-state-host2.service
After=mapper-wait@-xyz-openbmc_project-state-host2.service
Wants=mapper-wait@-xyz-openbmc_project-state-host3.service
After=mapper-wait@-xyz-openbmc_project-state-host3.service
Wants=mapper-wait@-xyz-openbmc_project-state-host4.service
After=mapper-wait@-xyz-openbmc_project-state-host4.service
Wants=mapper-wait@-xyz-openbmc_project-state-host5.service
After=mapper-wait@-xyz-openbmc_project-state-host5.service
Wants=mapper-wait@-xyz-openbmc_project-state-host6.service
After=mapper-wait@-xyz-openbmc_project-state-host6.service
Wants=mapper-wait@-xyz-openbmc_project-state-host6.service
After=mapper-wait@-xyz-openbmc_project-state-host6.service
Wants=mapper-wait@-xyz-openbmc_project-control-host1-acpi_power_state.service
After=mapper-wait@-xyz-openbmc_project-control-host1-acpi_power_state.service
Wants=mapper-wait@-xyz-openbmc_project-control-host2-acpi_power_state.service
After=mapper-wait@-xyz-openbmc_project-control-host2-acpi_power_state.service
Wants=mapper-wait@-xyz-openbmc_project-control-host3-acpi_power_state.service
After=mapper-wait@-xyz-openbmc_project-control-host3-acpi_power_state.service
Wants=mapper-wait@-xyz-openbmc_project-control-host4-acpi_power_state.service
After=mapper-wait@-xyz-openbmc_project-control-host4-acpi_power_state.service
Wants=mapper-wait@-xyz-openbmc_project-control-host5-acpi_power_state.service
After=mapper-wait@-xyz-openbmc_project-control-host5-acpi_power_state.service
Wants=mapper-wait@-xyz-openbmc_project-control-host6-acpi_power_state.service
After=mapper-wait@-xyz-openbmc_project-control-host6-acpi_power_state.service

[Service]
ExecStart=/usr/libexec/bletchley-host-state-monitor
SyslogIdentifier=bletchley-host-state-monitor
Type=simple
Restart=on-failure
RestartSec=1

[Install]
WantedBy=multi-user.target
