# meta-quanta/meta-gbs/recipes-gbs/hotswap-power-cycle/files/hotswap-power-cycle.service
[Unit]
Description=Power Cycle by Hotswap Controller

[Service]
EnvironmentFile=/run/psu_timedelay
ExecStartPre=sleep $PSU_HARDRESET_DELAY
ExecStart=gpioset gpiochip6 26=1
Restart=no
Type=oneshot
RemainAfterExit=true
StandardOutput=syslog

[Install]
WantedBy=gbmc-psu-hardreset.target
