# meta-quanta/meta-gbs/recipes-gbs/hotswap-power-cycle/files/hotswap-power-cycle.service
[Unit]
Description=Power Cycle by Hotswap Controller
After=gbmc-psu-hardreset.target

[Service]
Type=oneshot
ExecStart=/usr/bin/tray_powercycle.sh

[Install]
WantedBy=gbmc-psu-hardreset.target
