# yocto-poky/meta/recipes-bsp/acpid/acpid/acpid.service
[Unit]
Description=ACPI Event Daemon
After=syslog.target

[Service]
Type=forking
ExecStart=@SBINDIR@/acpid

[Install]
WantedBy=multi-user.target
