# meta-facebook/recipes-x86/chassis/x86-power-control/tiogapass/chassis-system-reset.service
[Unit]
Description=System unit to hard reset or system reset chassis

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=i2cset -y 7 0x45 0xd9 c

[Install]
WantedBy=chassis-system-reset.target
