# meta-amd/meta-ethanolx/recipes-amd/amd-fpga/files/ip-to-fpga.service
[Unit]
Description=Transfer IP address to the FPGA
BindsTo=sys-subsystem-net-devices-eth0.device
After=sys-subsystem-net-devices-eth0.device

[Service]
Type=simple
ExecStart=/usr/bin/ip-to-fpga.sh

[Install]
WantedBy=multi-user.target
