# meta-openbmc-machines/meta-openpower/common/recipes-phosphor/debug/openpower-debug-collector/op_checkstop_app.service
[Unit]
Description=Open power debug data collector

[Service]
ExecStart={sbindir}/checkstop_app
Type=oneshot
Restart=no
