# meta-openbmc-machines/meta-openpower/common/recipes-phosphor/debug/openpower-debug-collector/openpower-debug-collector-checkstop@.service
[Unit]
Description=OpenPOWER debug data collector for host checkstop
Before=obmc-host-quiesce@%i.target

[Service]
# Adding 5 seconds to allow OCC to collect debug data
# before processing checkstop
ExecStart={base_bindir}/sleep 5
ExecStart=/usr/bin/env checkstop_app
SyslogIdentifier=checkstop_app
Type=oneshot
Restart=no
