# meta-openbmc-machines/meta-openpower/common/recipes-phosphor/debug/openpower-debug-collector/openpower-debug-collector@.service
[Unit]
Description=OpenPOWER debug data collector
Before=obmc-quiesce-host@%i.target
Conflicts=obmc-host-quiesce@%i.target

[Service]
ExecStart={sbindir}/checkstop_app
Type=oneshot
Restart=no
