# meta-fii/meta-mori/recipes-mori/mori-fw-utility/mori-fw/mori-fw-ver.service
[Unit]
Description=Firmware Version Daemon

[Service]
RemainAfterExit=yes
Type=oneshot
ExecStart=/usr/libexec/mori-fw/mori-fw-ver.sh

[Install]
WantedBy=multi-user.target
