# meta-phosphor/recipes-phosphor/fans/obmc-control-fan/org.openbmc.control.Fans.service
[Unit]
Description=Phosphor Fan Control
Wants=obmc-fans-ready.target
After=obmc-fans-ready.target
Wants=obmc-fan-control.target
Before=obmc-fan-control.target

[Service]
Restart=always
ExecStart=/usr/bin/env fan_control.exe -f 6
SyslogIdentifier=fan_control.exe
Type=dbus
BusName={BUSNAME}
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
