# meta-phosphor/recipes-phosphor/fans/phosphor-fan/phosphor-fan-monitor-init@.service
[Unit]
Description=Phosphor Fan Monitor Initialization
Wants=obmc-power-on@%i.target
After=obmc-power-on@%i.target
Conflicts=obmc-chassis-powered-off@%i.target

[Service]
Restart=on-failure
ExecStart=/usr/bin/env phosphor-fan-monitor --init
SyslogIdentifier=phosphor-fan-monitor

[Install]
WantedBy=obmc-chassis-poweron@%i.target
