# meta-phosphor/recipes-phosphor/fans/phosphor-fan/phosphor-fan-monitor@.service
[Unit]
Description=Phosphor Fan Monitor Daemon
Conflicts=obmc-chassis-powered-off@%i.target

[Service]
Restart=on-failure
ExecStart=/usr/bin/env phosphor-fan-monitor --monitor
SyslogIdentifier=phosphor-fan-monitor

[Install]
RequiredBy=obmc-fan-control-ready@%i.target
