# meta-openbmc-machines/meta-openpower/meta-ingrasys/meta-zaius/recipes-phosphor/chassis/set-vio/set-vio.service
[Unit]
Description=Set VIO rail voltages
Before={SYSTEMD_DEFAULT_TARGET}

[Service]
ExecStart=/usr/sbin/i2cset -y 7 0x43 0x21 287 w
ExecStart=/usr/sbin/i2cset -y 8 0x42 0x21 287 w
Type=oneshot

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
