# meta-phosphor/common/recipes-phosphor/chassis/files/obmc-phosphor-chassisd.service
[Unit]
Description=Phosphor Chassis Control

[Service]
Restart=always
ExecStart={sbindir}/obmc-phosphor-chassisd

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
