# meta-ampere/meta-common/recipes-ac01/power-capping/ampere-power-cap/xyz.openbmc_project.AmpPowerCap.service
[Unit]
Description=Ampere Power Capping

[Service]
Restart=always
ExecStart=/usr/bin/ampere-power-cap
Type=simple

[Install]
WantedBy=multi-user.target
