# meta-phosphor/recipes-phosphor/dbus/phosphor-dbus-monitor/phosphor-dbus-monitor.service
[Unit]
Description=Phosphor DBus Monitor
After=obmc-mapper.target

[Service]
Restart=always
ExecStart=/usr/bin/env phosphor-dbus-monitor
SyslogIdentifier=phosphor-dbus-monitor

[Install]
WantedBy=multi-user.target
