# meta-phosphor/recipes-phosphor/dbus/phosphor-dbus-monitor/phosphor-msl-verify.service
[Unit]
Description=Phosphor MSL Verification
After=mapper-wait@-xyz-openbmc_project-inventory.service
Wants=mapper-wait@-xyz-openbmc_project-inventory.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/env phosphor-msl-verify
SyslogIdentifier=phosphor-msl-verify
