# meta-phosphor/common/recipes-phosphor/chassis/obmc-button-power/org.openbmc.buttons.Power@.service
[Unit]
Description=Phosphor Power Button%i
Wants=mapper-wait@-org-openbmc-managers-System.service
After=mapper-wait@-org-openbmc-managers-System.service

[Service]
Restart=always
ExecStart=/usr/bin/env button_power.exe
SyslogIdentifier=button_power.exe
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy={SYSTEMD_DEFAULT_TARGET}
